# ts/5.1/packages/xfwm4/build/extra/etc/xdg/autostart/wait-for-tray.desktop
[Desktop Entry]
Name=wait-for-tray
Comment=wait-for-tray
Exec=/bin/wait-for-tray
Terminal=false
Type=Application
Categories=
NoDisplay=true
OnlyShowIn=GNOME;XFCE;Unity;

