# ts/5.1/packages/xfwm4/build/extra/etc/xdg/autostart/xfwm4-thinstation.desktop
[Desktop Entry]
Version=1.0
Name=Xfsettings-thinstation
Exec=/bin/xfwm4-thinstation
Icon=applications-system
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=XFCE;
