# ts/build/packages/plymouth/build/extra/lib/systemd/system/plymouth-quit.service
[Unit]
Description=Terminate Plymouth Boot Screen
After=rc-local.service plymouth-start.service systemd-user-sessions.service session-setup.service

[Service]
ExecStart=-/usr/bin/plymouth quit
Type=oneshot
TimeoutSec=20
