# ts/ports/contrib/ntp/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/ntp/
-rw-r--r--	root/root	etc/ntp/ntp.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/ntpd
-rwxr-xr-x	root/root	etc/rc.d/ntpdate
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/calc_tickadj
-rwxr-xr-x	root/root	usr/bin/ntp-keygen
-rwxr-xr-x	root/root	usr/bin/ntp-wait
-rwxr-xr-x	root/root	usr/bin/ntpd
-rwxr-xr-x	root/root	usr/bin/ntpdate
-rwxr-xr-x	root/root	usr/bin/ntpdc
-rwxr-xr-x	root/root	usr/bin/ntpq
-rwxr-xr-x	root/root	usr/bin/ntptime
-rwxr-xr-x	root/root	usr/bin/ntptrace
-rwxr-xr-x	root/root	usr/bin/sntp
-rwxr-xr-x	root/root	usr/bin/tickadj
-rwxr-xr-x	root/root	usr/bin/update-leap
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/calc_tickadj.1.gz
-rw-r--r--	root/root	usr/man/man1/ntp-keygen.1.gz
-rw-r--r--	root/root	usr/man/man1/ntp-wait.1.gz
-rw-r--r--	root/root	usr/man/man1/ntpd.1.gz
-rw-r--r--	root/root	usr/man/man1/ntpdc.1.gz
-rw-r--r--	root/root	usr/man/man1/ntpq.1.gz
-rw-r--r--	root/root	usr/man/man1/ntptrace.1.gz
-rw-r--r--	root/root	usr/man/man1/sntp.1.gz
-rw-r--r--	root/root	usr/man/man1/update-leap.1.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/ntp.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/ntp.keys.5.gz
drwxr-xr-x	root/root	usr/man/man8/
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/ntp/
drwxr-xr-x	root/root	usr/share/ntp/lib/
drwxr-xr-x	root/root	usr/share/ntp/lib/NTP/
-rw-r--r--	root/root	usr/share/ntp/lib/NTP/Util.pm
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/ntp/
-rw-r--r--	root/root	var/lib/ntp/ntp.drift
drwxr-xr-x	root/root	var/lib/ntp/stats/
drwxr-xr-x	root/root	var/log/
-rw-r--r--	root/root	var/log/ntp.log (EMPTY)
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	root/root	var/run/ntp/
