# ts/ports/binary-core/firmware-ap6214a/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/firmware/
drwxr-xr-x	root/root	etc/firmware/ampak-ap6214a/
-rw-r--r--	root/root	etc/firmware/ampak-ap6214a/BCM43430A1.hcd
drwxr-xr-x	root/root	etc/init.d/
-rw-r--r--	root/root	etc/init.d/ap6214a
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/firmware/
drwxr-xr-x	root/root	lib/firmware/ampak-ap6214a/
drwxr-xr-x	root/root	lib/firmware/ampak-ap6214a/brcm/
-rw-r--r--	root/root	lib/firmware/ampak-ap6214a/brcm/brcmfmac43430-sdio.bin
-rw-r--r--	root/root	lib/firmware/ampak-ap6214a/brcm/brcmfmac43430-sdio.txt
drwxr-xr-x	root/root	lib/modprobe.d/
-rw-r--r--	root/root	lib/modprobe.d/brcmfmac-ampak.conf
