# ts/ports/opt/cups-filters/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/cups/
-rw-r--r--	root/root	etc/cups/cups-browsed.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/cups-browsed
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
lrwxrwxrwx	root/root	usr/bin/driverless -> ../lib/cups/driver/driverless
lrwxrwxrwx	root/root	usr/bin/driverless-fax -> ../lib/cups/driver/driverless-fax
lrwxrwxrwx	root/root	usr/bin/foomatic-rip -> ../lib/cups/filter/foomatic-rip
-rwxr-xr-x	root/root	usr/bin/ttfread
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/cupsfilters/
-rw-r--r--	root/root	usr/include/cupsfilters/colord.h
-rw-r--r--	root/root	usr/include/cupsfilters/colormanager.h
-rw-r--r--	root/root	usr/include/cupsfilters/driver.h
-rw-r--r--	root/root	usr/include/cupsfilters/image.h
-rw-r--r--	root/root	usr/include/cupsfilters/ipp.h
-rw-r--r--	root/root	usr/include/cupsfilters/pdftoippprinter.h
-rw-r--r--	root/root	usr/include/cupsfilters/ppdgenerator.h
-rw-r--r--	root/root	usr/include/cupsfilters/raster.h
drwxr-xr-x	root/root	usr/include/fontembed/
-rw-r--r--	root/root	usr/include/fontembed/bitset.h
-rw-r--r--	root/root	usr/include/fontembed/embed.h
-rw-r--r--	root/root	usr/include/fontembed/fontfile.h
-rw-r--r--	root/root	usr/include/fontembed/iofn.h
-rw-r--r--	root/root	usr/include/fontembed/sfnt.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/backend/
-rwxr-xr-x	root/root	usr/lib/cups/backend/beh
lrwxrwxrwx	root/root	usr/lib/cups/backend/driverless -> ../driver/driverless
lrwxrwxrwx	root/root	usr/lib/cups/backend/driverless-fax -> ../driver/driverless-fax
-rwx------	root/root	usr/lib/cups/backend/implicitclass
-rwxr-xr-x	root/root	usr/lib/cups/backend/parallel
-rwxr-xr-x	root/root	usr/lib/cups/backend/serial
drwxr-xr-x	root/root	usr/lib/cups/driver/
-rwxr-xr-x	root/root	usr/lib/cups/driver/driverless
-rwxr-xr-x	root/root	usr/lib/cups/driver/driverless-fax
drwxr-xr-x	root/root	usr/lib/cups/filter/
-rwxr-xr-x	root/root	usr/lib/cups/filter/bannertopdf
-rwxr-xr-x	root/root	usr/lib/cups/filter/commandtoescpx
-rwxr-xr-x	root/root	usr/lib/cups/filter/commandtopclx
-rwxr-xr-x	root/root	usr/lib/cups/filter/foomatic-rip
-rwxr-xr-x	root/root	usr/lib/cups/filter/gstopdf
-rwxr-xr-x	root/root	usr/lib/cups/filter/gstopxl
-rwxr-xr-x	root/root	usr/lib/cups/filter/gstoraster
-rwxr-xr-x	root/root	usr/lib/cups/filter/imagetopdf
-rwxr-xr-x	root/root	usr/lib/cups/filter/imagetops
-rwxr-xr-x	root/root	usr/lib/cups/filter/imagetoraster
-rwxr-xr-x	root/root	usr/lib/cups/filter/mupdftoraster
-rwxr-xr-x	root/root	usr/lib/cups/filter/pdftopdf
-rwxr-xr-x	root/root	usr/lib/cups/filter/pdftops
-rwxr-xr-x	root/root	usr/lib/cups/filter/pdftoraster
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertoescpx
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertopclm
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertopclx
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertopdf
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertops
-rwxr-xr-x	root/root	usr/lib/cups/filter/sys5ippprinter
-rwxr-xr-x	root/root	usr/lib/cups/filter/texttopdf
-rwxr-xr-x	root/root	usr/lib/cups/filter/texttops
-rwxr-xr-x	root/root	usr/lib/cups/filter/texttotext
-rw-r--r--	root/root	usr/lib/libcupsfilters.a
-rwxr-xr-x	root/root	usr/lib/libcupsfilters.la
lrwxrwxrwx	root/root	usr/lib/libcupsfilters.so -> libcupsfilters.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libcupsfilters.so.1 -> libcupsfilters.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libcupsfilters.so.1.0.0
-rw-r--r--	root/root	usr/lib/libfontembed.a
-rwxr-xr-x	root/root	usr/lib/libfontembed.la
lrwxrwxrwx	root/root	usr/lib/libfontembed.so -> libfontembed.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libfontembed.so.1 -> libfontembed.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libfontembed.so.1.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libcupsfilters.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libfontembed.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/cups-browsed
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/cups/
drwxr-xr-x	root/root	usr/share/cups/banners/
-rw-r--r--	root/root	usr/share/cups/banners/classified
-rw-r--r--	root/root	usr/share/cups/banners/confidential
-rw-r--r--	root/root	usr/share/cups/banners/form
-rw-r--r--	root/root	usr/share/cups/banners/secret
-rw-r--r--	root/root	usr/share/cups/banners/standard
-rw-r--r--	root/root	usr/share/cups/banners/topsecret
-rw-r--r--	root/root	usr/share/cups/banners/unclassified
drwxr-xr-x	root/root	usr/share/cups/charsets/
lrwxrwxrwx	root/root	usr/share/cups/charsets/pdf.utf-8 -> pdf.utf-8.simple
-rw-r--r--	root/root	usr/share/cups/charsets/pdf.utf-8.heavy
-rw-r--r--	root/root	usr/share/cups/charsets/pdf.utf-8.simple
drwxr-xr-x	root/root	usr/share/cups/data/
-rw-r--r--	root/root	usr/share/cups/data/classified.pdf
-rw-r--r--	root/root	usr/share/cups/data/confidential.pdf
-rw-r--r--	root/root	usr/share/cups/data/default-testpage.pdf
-rw-r--r--	root/root	usr/share/cups/data/default.pdf
-rw-r--r--	root/root	usr/share/cups/data/form_english.pdf
-rw-r--r--	root/root	usr/share/cups/data/form_english_in.odt
-rw-r--r--	root/root	usr/share/cups/data/form_russian.pdf
-rw-r--r--	root/root	usr/share/cups/data/form_russian_in.odt
-rw-r--r--	root/root	usr/share/cups/data/secret.pdf
-rw-r--r--	root/root	usr/share/cups/data/standard.pdf
-rw-r--r--	root/root	usr/share/cups/data/testprint
-rw-r--r--	root/root	usr/share/cups/data/topsecret.pdf
-rw-r--r--	root/root	usr/share/cups/data/unclassified.pdf
drwxr-xr-x	root/root	usr/share/cups/drv/
-rw-r--r--	root/root	usr/share/cups/drv/cupsfilters.drv
drwxr-xr-x	root/root	usr/share/cups/mime/
-rw-r--r--	root/root	usr/share/cups/mime/cupsfilters-ghostscript.convs
-rw-r--r--	root/root	usr/share/cups/mime/cupsfilters-mupdf.convs
-rw-r--r--	root/root	usr/share/cups/mime/cupsfilters-poppler.convs
-rw-r--r--	root/root	usr/share/cups/mime/cupsfilters.convs
-rw-r--r--	root/root	usr/share/cups/mime/cupsfilters.types
drwxr-xr-x	root/root	usr/share/cups/ppdc/
-rw-r--r--	root/root	usr/share/cups/ppdc/escp.h
-rw-r--r--	root/root	usr/share/cups/ppdc/pcl.h
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/driverless.1.gz
-rw-r--r--	root/root	usr/share/man/man1/foomatic-rip.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/cups-browsed.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/cups-browsed.8.gz
drwxr-xr-x	root/root	usr/share/ppd/
drwxr-xr-x	root/root	usr/share/ppd/cupsfilters/
-rw-r--r--	root/root	usr/share/ppd/cupsfilters/Fuji_Xerox-DocuPrint_CM305_df-PDF.ppd
-rw-r--r--	root/root	usr/share/ppd/cupsfilters/Generic-PDF_Printer-PDF.ppd
-rw-r--r--	root/root	usr/share/ppd/cupsfilters/HP-Color_LaserJet_CM3530_MFP-PDF.ppd
-rw-r--r--	root/root	usr/share/ppd/cupsfilters/Ricoh-PDF_Printer-PDF.ppd
-rw-r--r--	root/root	usr/share/ppd/cupsfilters/pxlcolor.ppd
-rw-r--r--	root/root	usr/share/ppd/cupsfilters/pxlmono.ppd
