# ts/build/packages/vmview/build/extra/etc/systemd/system/ftscanhv.service
# Automatically generated by systemd-sysv-generator

[Unit]
Description=VMware Horizon Client Scanner Service
Documentation=man:systemd-sysv-generator(8)
SourcePath=/etc/init.d/ftscanhv
Requires=vmware.usbarbitrator.service
After=vmware.usbarbitrator.service pkg.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/ftscanhv start
ExecStop=/etc/init.d/ftscanhv stop
