# ts/build/packages/vmview/build/extra/etc/systemd/system/ftsprhv.service
# Automatically generated by systemd-sysv-generator

[Unit]
Description=VMware Horizon Client USB Sharing Service
Documentation=man:systemd-sysv-generator(8)
SourcePath=/etc/init.d/ftsprhv
Requires=vmware.usbarbitrator.service
After=vmware.usbarbitrator.service pkg.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/ftsprhv start
ExecStop=/etc/init.d/ftsprhv stop
