# ts/ports/core/exim/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/exim/
-rw-r--r--	root/root	etc/exim/aliases
-rw-r--r--	root/root	etc/exim/exim.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/exim
drwxr-xr-x	root/root	etc/ssl/
drwxr-xr-x	root/root	etc/ssl/certs/
-rw-------	root/root	etc/ssl/certs/exim.crt (EMPTY)
drwxr-xr-x	root/root	etc/ssl/keys/
-rw-------	root/root	etc/ssl/keys/exim.key (EMPTY)
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
lrwxrwxrwx	root/root	usr/bin/mailq -> ../sbin/exim
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/exicyclog
-rwxr-xr-x	root/root	usr/sbin/exigrep
-rwsr-xr-x	root/root	usr/sbin/exim
-rwxr-xr-x	root/root	usr/sbin/exim_checkaccess
-rwxr-xr-x	root/root	usr/sbin/exim_dbmbuild
-rwxr-xr-x	root/root	usr/sbin/exim_dumpdb
-rwxr-xr-x	root/root	usr/sbin/exim_fixdb
-rwxr-xr-x	root/root	usr/sbin/exim_lock
-rwxr-xr-x	root/root	usr/sbin/exim_tidydb
-rwxr-xr-x	root/root	usr/sbin/eximstats
-rwxr-xr-x	root/root	usr/sbin/exinext
-rwxr-xr-x	root/root	usr/sbin/exipick
-rwxr-xr-x	root/root	usr/sbin/exiqgrep
-rwxr-xr-x	root/root	usr/sbin/exiqsumm
-rwxr-xr-x	root/root	usr/sbin/exiwhat
lrwxrwxrwx	root/root	usr/sbin/sendmail -> exim
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/exim.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-x---	mail/mail	var/log/exim/
-rw-r-----	mail/mail	var/log/exim/exim_mainlog (EMPTY)
-rw-r-----	mail/mail	var/log/exim/exim_paniclog (EMPTY)
-rw-r-----	mail/mail	var/log/exim/exim_rejectlog (EMPTY)
drwxr-xr-x	root/root	var/spool/
drwxr-x---	mail/mail	var/spool/exim/
