# ts/ports/components/easytether/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/easytether/
drwxr-xr-x	root/root	etc/network/
drwxr-xr-x	root/root	etc/network/interfaces.d/
-rw-r--r--	root/root	etc/network/interfaces.d/tun-easytether
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/systemd/
drwxr-xr-x	root/root	lib/systemd/network/
-rw-r--r--	root/root	lib/systemd/network/99-tun-easytether.network
drwxr-xr-x	root/root	lib/systemd/system/
-rw-r--r--	root/root	lib/systemd/system/easytether-usb@.service
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/99-easytether-usb.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/easytether-bluetooth
-rwxr-xr-x	root/root	usr/bin/easytether-local
-rwxr-xr-x	root/root	usr/bin/easytether-usb
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/easytether/
-rw-r--r--	root/root	usr/share/doc/easytether/changelog.gz
-rw-r--r--	root/root	usr/share/doc/easytether/copyright
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/easytether/
drwxr-xr-x	root/root	var/lib/easytether/empty/
