# ts/ports/xorg-opt/modemmanager/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.ModemManager1.conf
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-cinterion-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-ericsson-mbm.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-huawei-net-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-longcheer-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-mtk-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-nokia-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-pcmcia-device-blacklist.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-platform-serial-whitelist.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-simtech-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-telit-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-usb-device-blacklist.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-usb-serial-adapters-greylist.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-x22x-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/77-mm-zte-port-types.rules
-rw-r--r--	root/root	lib/udev/rules.d/80-mm-candidate.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/mmcli
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/ModemManager/
-rw-r--r--	root/root	usr/include/ModemManager/ModemManager-enums.h
-rw-r--r--	root/root	usr/include/ModemManager/ModemManager-errors.h
-rw-r--r--	root/root	usr/include/ModemManager/ModemManager-names.h
-rw-r--r--	root/root	usr/include/ModemManager/ModemManager-version.h
-rw-r--r--	root/root	usr/include/ModemManager/ModemManager.h
drwxr-xr-x	root/root	usr/include/libmm-glib/
-rw-r--r--	root/root	usr/include/libmm-glib/libmm-glib.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-bearer-ip-config.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-bearer-properties.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-bearer.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-cdma-manual-activation-properties.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-enums-types.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-errors-types.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-firmware-properties.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-gdbus-bearer.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-gdbus-manager.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-gdbus-modem.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-gdbus-sim.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-gdbus-sms.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-helper-types.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-location-3gpp.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-location-cdma-bs.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-location-common.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-location-gps-nmea.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-location-gps-raw.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-manager.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-3gpp-ussd.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-3gpp.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-cdma.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-firmware.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-location.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-messaging.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-oma.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-signal.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-simple.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem-time.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-modem.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-network-timezone.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-object.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-signal.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-sim.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-simple-connect-properties.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-simple-status.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-sms-properties.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-sms.h
-rw-r--r--	root/root	usr/include/libmm-glib/mm-unlock-retries.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/ModemManager/
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-altair-lte.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-altair-lte.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-altair-lte.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-anydata.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-anydata.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-anydata.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-cinterion.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-cinterion.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-cinterion.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-generic.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-generic.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-generic.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-gobi.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-gobi.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-gobi.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-hso.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-hso.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-hso.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-huawei.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-huawei.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-huawei.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-iridium.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-iridium.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-iridium.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-linktop.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-linktop.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-linktop.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-longcheer.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-longcheer.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-longcheer.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-mbm.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-mbm.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-mbm.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-motorola.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-motorola.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-motorola.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-mtk.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-mtk.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-mtk.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-nokia-icera.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-nokia-icera.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-nokia-icera.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-nokia.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-nokia.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-nokia.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-novatel-lte.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-novatel-lte.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-novatel-lte.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-novatel.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-novatel.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-novatel.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-option.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-option.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-option.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-pantech.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-pantech.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-pantech.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-samsung.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-samsung.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-samsung.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-sierra.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-sierra.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-sierra.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-simtech.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-simtech.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-simtech.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-telit.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-telit.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-telit.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-via.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-via.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-via.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-wavecom.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-wavecom.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-wavecom.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-x22x.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-x22x.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-x22x.so
-rw-r--r--	root/root	usr/lib/ModemManager/libmm-plugin-zte.a
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-zte.la
-rwxr-xr-x	root/root	usr/lib/ModemManager/libmm-plugin-zte.so
-rw-r--r--	root/root	usr/lib/libmm-glib.a
-rwxr-xr-x	root/root	usr/lib/libmm-glib.la
lrwxrwxrwx	root/root	usr/lib/libmm-glib.so -> libmm-glib.so.0.2.0
lrwxrwxrwx	root/root	usr/lib/libmm-glib.so.0 -> libmm-glib.so.0.2.0
-rwxr-xr-x	root/root	usr/lib/libmm-glib.so.0.2.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/ModemManager.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/mm-glib.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/ModemManager.service
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/ModemManager
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/interfaces/
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Bearer.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Firmware.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Location.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Messaging.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Modem3gpp.Ussd.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Modem3gpp.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.ModemCdma.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Oma.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Signal.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Simple.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Time.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Sim.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Sms.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/wip-org.freedesktop.ModemManager1.Modem.Contacts.xml
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.ModemManager1.service
drwxr-xr-x	root/root	usr/share/icons/
drwxr-xr-x	root/root	usr/share/icons/hicolor/
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/22x22/apps/ModemManager.png
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/ModemManager.8.gz
-rw-r--r--	root/root	usr/share/man/man8/mmcli.8.gz
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/org.freedesktop.ModemManager1.policy
