# ts/build/packages/sshd/etc/systemd/system/dropbear.service
[Unit]
Description=Dropbear SSH server
After=network.target pkg.service

[Service]
Type=forking
ExecStart=/bin/dropbear -r /etc/dropbear/dropbear_rsa_host_key
PIDFile=/run/dropbear.pid

[Install]
WantedBy=multi-user.target
