# ts/build/packages/telnetd/etc/systemd/system/telnetd-init.service
[Unit]
Description=ThinStation Telnet Setup
Wants=network-online.target
After=network-files.service
Before=shutdown.target xinetd.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/telnetd-init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target

