# ts/build/packages/lpr/etc/systemd/system/lpr.service
[Unit]
Description=Thistation lpr
After=profile-setup.service pkg.service
Before=network-pre.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/lpr init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
