# ts/build/packages/ntpd/build/extra/etc/systemd/system/ntpd.service.d/override.conf
[Unit]
After=fastboot.service NetworkManager-wait-online.service
