# ts/build/packages/ntpd/build/extra/etc/systemd/system/ntpd.service.d/override.conf
[Unit]
After=fastboot.service NetworkManager-wait-online.service

[Service]
EnvironmentFile=/etc/thinstation.env
ExecStartPre=/etc/init.d/ntpconf init
ExecStartPost=/etc/init.d/ntpconf update_rtc
ExecStopPost=/etc/init.d/ntpconf update_rtc

