# ts/5.1/packages/devstation/etc/xinetd.d/http
service http
{
	flags			= REUSE
	socket_type		= stream        
	wait			= no
	user			= root
	server			= /sbin/httpd
	server_args		= -i -h /var/lib/tftpboot
        log_type        	= FILE /var/log/httpd
	log_on_failure		+= USERID
	disable			= no
}
