# ts/ports/components/dmraid/.footprint
drwxr-xr-x	root/root	sbin/
-rwxr-xr-x	root/root	sbin/dmevent_tool
-rwxr-xr-x	root/root	sbin/dmraid
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/dmraid/
-rw-r--r--	root/root	usr/include/dmraid/display.h
-rw-r--r--	root/root	usr/include/dmraid/dmraid.h
-rw-r--r--	root/root	usr/include/dmraid/dmreg.h
-rw-r--r--	root/root	usr/include/dmraid/format.h
-rw-r--r--	root/root	usr/include/dmraid/lib_context.h
-rw-r--r--	root/root	usr/include/dmraid/list.h
-rw-r--r--	root/root	usr/include/dmraid/locking.h
-rw-r--r--	root/root	usr/include/dmraid/metadata.h
-rw-r--r--	root/root	usr/include/dmraid/misc.h
-rw-r--r--	root/root	usr/include/dmraid/reconfig.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/device-mapper/
-rwxr-xr-x	root/root	usr/lib/device-mapper/libdmraid-events-isw.so
drwxr-xr-x	root/root	usr/lib/initcpio/
drwxr-xr-x	root/root	usr/lib/initcpio/hooks/
-rw-r--r--	root/root	usr/lib/initcpio/hooks/dmraid
drwxr-xr-x	root/root	usr/lib/initcpio/install/
-rw-r--r--	root/root	usr/lib/initcpio/install/dmraid
lrwxrwxrwx	root/root	usr/lib/libdmraid-events-isw.so -> device-mapper/libdmraid-events-isw.so
-rw-r--r--	root/root	usr/lib/libdmraid.a
lrwxrwxrwx	root/root	usr/lib/libdmraid.so -> libdmraid.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libdmraid.so.1 -> libdmraid.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libdmraid.so.1.0.0
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/dmraid.service
drwxr-xr-x	root/root	usr/lib/tmpfiles.d/
-rw-r--r--	root/root	usr/lib/tmpfiles.d/dmraid.conf
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/dmevent_tool.8.gz
-rw-r--r--	root/root	usr/share/man/man8/dmraid.8.gz
