# ts/ports/core/dcron/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/cron.d/
drwxr-xr-x	root/root	etc/cron/
drwxr-xr-x	root/root	etc/cron/daily/
drwxr-xr-x	root/root	etc/cron/hourly/
drwxr-xr-x	root/root	etc/cron/monthly/
drwxr-xr-x	root/root	etc/cron/weekly/
lrwxrwxrwx	root/root	etc/crontab -> ../var/spool/cron/crontabs/root
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/crond
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwsr-x---	root/users	usr/bin/crontab
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/crontab.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/crond.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwx------	root/root	usr/sbin/crond
-rwxr-xr-x	root/root	usr/sbin/runjobs
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/spool/
drwxr-xr-x	root/root	var/spool/cron/
drwxr-xr-x	root/root	var/spool/cron/cronstamps/
drwxr-xr-x	root/root	var/spool/cron/crontabs/
-rw-------	root/root	var/spool/cron/crontabs/root
