# ts/5.1/packages/base/etc/inittab
# /etc/inittab init(8) configuration for BusyBox

tty1::sysinit:-/etc/thinstation.init

tty1::respawn:-/bin/start-session 0
tty2::askfirst:/bin/getty tty2 9600
tty3::askfirst:/bin/getty tty3 9600

tty1::shutdown:/bin/shutdown

::ctrlaltdel:/bin/reboot
