# ts/build/packages/sound-nasd/etc/systemd/system/sound-nasd.service
[Unit]
Description=ThinStation sound-nasd
After=profile-setup.service pkg.service
Before=display-manager.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/sound-nasd init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
