# ts/ports/components/dnscrypt-proxy/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dnscrypt-proxy/
-rw-r--r--	root/root	etc/dnscrypt-proxy/blacklist.txt
-rw-r--r--	root/root	etc/dnscrypt-proxy/cloaking-rules.txt
-rw-r--r--	root/root	etc/dnscrypt-proxy/dnscrypt-proxy.toml
-rw-r--r--	root/root	etc/dnscrypt-proxy/forwarding-rules.txt
-rw-r--r--	root/root	etc/dnscrypt-proxy/ip-blacklist.txt (EMPTY)
-rw-r--r--	root/root	etc/dnscrypt-proxy/whitelist.txt
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/dnscrypt-proxy
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/dnscrypt-proxy.service
-rw-r--r--	root/root	usr/lib/systemd/system/dnscrypt-proxy.socket
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/dnscrypt-proxy/
-rw-r--r--	root/root	usr/share/doc/dnscrypt-proxy/ChangeLog
-rw-r--r--	root/root	usr/share/doc/dnscrypt-proxy/README.md
drwxr-xr-x	root/root	usr/share/licenses/
drwxr-xr-x	root/root	usr/share/licenses/dnscrypt-proxy/
-rw-r--r--	root/root	usr/share/licenses/dnscrypt-proxy/LICENSE
