# ts/ports/components/hostapd/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/hostapd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/hostapd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/hostapd
-rwxr-xr-x	root/root	usr/sbin/hostapd_cli
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/hostapd_cli.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/hostapd.8.gz
