# ts/build/packages/networkmanager/build/extra/etc/NetworkManager/dispatcher.d/00-update-status
#!/bin/sh

. /etc/thinstation.global

INTERFACE=$1
ACTION=$2

if [ "$INTERFACE" != "lo" ]; then

	case $ACTION in
	pre-up)
		CLIENT_MAC=`cat /sys/class/net/$INTERFACE/address | sed 's/://g'`
		CLIENT_NAME=$HOSTNAME
	;;
	up|dhcp4-change)
		if [ -e /var/log/net/$INTERFACE ]; then
			. /var/log/net/$INTERFACE
		fi
		if [ -n "$DHCP4_HOST_NAME" ]; then
			CLIENT_NAME="$DHCP4_HOST_NAME"
			hostnamectl set-hostname $CLIENT_NAME
			sed -i -e '/localhost/d' /etc/hosts
                        echo "127.0.0.1 `hostname` localhost" >> /etc/hosts

			# Update the .Xauthority with the new hostname
			if `which xauth >/dev/null` && [ -n "`xauth list`" ]; then
				# Extract the old key
				OLDHOST=`xauth list | awk '{print $1}' | head -1`
				KEY=`xauth list | awk '{print $3}' | head -1`

				# Get the current hostname
				NEWHOST=`hostname`/unix:0

				# Add new entry
				xauth add $NEWHOST MIT-MAGIC-COOKIE-1 $KEY

				# Remove the old entry
				xauth remove "$OLDHOST"
			fi
		fi
		NETWORKUP=TRUE
		gw=`echo $IP4_ADDRESS_0 | cut -d/ -f2 | cut -d' ' -f2`
		ip=`echo $IP4_ADDRESS_0 | cut -d/ -f1`

		sed -i -e "/NET_DEVICE=/d" $TS_RUNTIME
		echo "NET_DEVICE=$INTERFACE" >> $TS_RUNTIME
	;;
	down)
		if [ -e /var/log/net/$INTERFACE ]; then
			. /var/log/net/$INTERFACE
		fi
		NETWORKUP=FALSE
	;;
	*)
		exit
	;;
	esac

	mkdir -p /var/log/net
	echo "DEVTYPE=\"$CONNECTION_ID\""			> /var/log/net/$INTERFACE
	echo "CLIENT_NAME=$CLIENT_NAME"			>> /var/log/net/$INTERFACE
	echo "CLIENT_MAC=$CLIENT_MAC"			>> /var/log/net/$INTERFACE
	echo "CLIENT_IP=$ip"				>> /var/log/net/$INTERFACE
	echo "CLIENT_GATEWAY=$gw"			>> /var/log/net/$INTERFACE
	echo "SUBNET=$DHCP4_NETWORK_NUMBER"		>> /var/log/net/$INTERFACE
	echo "SERVER_IP=$DHCP4_NEXT_SERVER"		>> /var/log/net/$INTERFACE
	echo "NETWORKUP=$NETWORKUP"			>> /var/log/net/$INTERFACE
	echo "NETMASK_SIZE=$DHCP4_SUBNET_MASK"		>> /var/log/net/$INTERFACE
	echo "DHCP_SERVER_ID=$DHCP4_DHCP_SERVER_IDENTIFIER">> /var/log/net/$INTERFACE
	echo "NET_USE=$NET_USE"				>> /var/log/net/$INTERFACE
	echo "NET_DHCP_TIMEOUT=$NET_DHCP_TIMEOUT"	>> /var/log/net/$INTERFACE
fi
