# ts/ports/components/haveged/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/haveged/
-rw-r--r--	root/root	usr/include/haveged/havege.h
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/libhavege.la
lrwxrwxrwx	root/root	usr/lib/libhavege.so -> libhavege.so.1.1.0
lrwxrwxrwx	root/root	usr/lib/libhavege.so.1 -> libhavege.so.1.1.0
-rwxr-xr-x	root/root	usr/lib/libhavege.so.1.1.0
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/haveged.service
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/haveged
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man3/
-rw-r--r--	root/root	usr/share/man/man3/libhavege.3.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/haveged.8.gz
