# ts/5.1/packages/iphone/lib/udev/rules.d/90-iphone-tether.rules
# udev rules for setting correct configuration and pairing on tethered iPhones
ATTR{idVendor}!="05ac", GOTO="ipheth_rules_end"

# Execute pairing program when appropriate
ACTION=="add", SUBSYSTEM=="net", ENV{ID_USB_DRIVER}=="ipheth", SYMLINK+="iphone", RUN+="/bin/ideviceinfo -k com.apple.mobile.tethered_sync"

LABEL="ipheth_rules_end"
