# ts/5.2/packages/networkmanager/build/extra/etc/udev/rules.d/99-tsnet.rules
SUBSYSTEM!="net", GOTO="net_end"
RUN+="/etc/udev/scripts/net.sh"
LABEL="net_end"
