# ts/build/packages/autonet/etc/systemd/system/autonet.service
[Unit]
Description=Configure Network Connections Automatically
DefaultDependencies=no
After=systemd-udev-settle.service
Before=sysinit.target choose-network.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/autonet
SyslogIdentifier=thinstation
