# ts/ports/opt/php-fcgi/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/php/
-rwxr-xr-x	root/root	etc/php/php.ini
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/php-fcgi
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/php-cgi
-rwxr-xr-x	root/root	usr/bin/php-config
-rwxr-xr-x	root/root	usr/bin/phpize
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/php/
drwxr-xr-x	root/root	usr/lib/php/extensions/
-rwxr-xr-x	root/root	usr/lib/php/extensions/opcache.so
