# ts/ports/contrib/spice/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/chezdav
-rwxr-xr-x	root/root	usr/bin/spicy
-rwxr-xr-x	root/root	usr/bin/spicy-screenshot
-rwxr-xr-x	root/root	usr/bin/spicy-stats
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/spice-client-glib-2.0/
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-cursor.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-display.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-inputs.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-main.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-playback.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-port.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-record.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-smartcard.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-usbredir.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/channel-webdav.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/smartcard-manager.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-audio.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-channel-enums.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-channel.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-client.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-glib-enums.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-option.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-session.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-types.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-uri.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-util.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/spice-version.h
-rw-r--r--	root/root	usr/include/spice-client-glib-2.0/usb-device-manager.h
drwxr-xr-x	root/root	usr/include/spice-client-gtk-3.0/
-rw-r--r--	root/root	usr/include/spice-client-gtk-3.0/spice-grabsequence.h
-rw-r--r--	root/root	usr/include/spice-client-gtk-3.0/spice-gtk-session.h
-rw-r--r--	root/root	usr/include/spice-client-gtk-3.0/spice-widget-enums.h
-rw-r--r--	root/root	usr/include/spice-client-gtk-3.0/spice-widget.h
-rw-r--r--	root/root	usr/include/spice-client-gtk-3.0/usb-device-widget.h
drwxr-xr-x	root/root	usr/include/spice-controller/
-rw-r--r--	root/root	usr/include/spice-controller/spice-controller.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/girepository-1.0/
-rw-r--r--	root/root	usr/lib/girepository-1.0/SpiceClientGLib-2.0.typelib
-rw-r--r--	root/root	usr/lib/girepository-1.0/SpiceClientGtk-3.0.typelib
-rwxr-xr-x	root/root	usr/lib/libspice-client-glib-2.0.la
lrwxrwxrwx	root/root	usr/lib/libspice-client-glib-2.0.so -> libspice-client-glib-2.0.so.8.5.0
lrwxrwxrwx	root/root	usr/lib/libspice-client-glib-2.0.so.8 -> libspice-client-glib-2.0.so.8.5.0
-rwxr-xr-x	root/root	usr/lib/libspice-client-glib-2.0.so.8.5.0
-rwxr-xr-x	root/root	usr/lib/libspice-client-gtk-3.0.la
lrwxrwxrwx	root/root	usr/lib/libspice-client-gtk-3.0.so -> libspice-client-gtk-3.0.so.4.0.0
lrwxrwxrwx	root/root	usr/lib/libspice-client-gtk-3.0.so.4 -> libspice-client-gtk-3.0.so.4.0.0
-rwxr-xr-x	root/root	usr/lib/libspice-client-gtk-3.0.so.4.0.0
-rwxr-xr-x	root/root	usr/lib/libspice-controller.la
lrwxrwxrwx	root/root	usr/lib/libspice-controller.so -> libspice-controller.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libspice-controller.so.0 -> libspice-controller.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libspice-controller.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/spice-client-glib-2.0.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/spice-client-gtk-3.0.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/spice-controller.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/spice-webdavd.service
-rw-r--r--	root/root	usr/lib/systemd/system/spice-webdavd.target
drwxr-xr-x	root/root	usr/lib/udev/
drwxr-xr-x	root/root	usr/lib/udev/rules.d/
-rw-r--r--	root/root	usr/lib/udev/rules.d/70-spice-webdavd.rules
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/chezdav.1.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/spice-webdavd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/gir-1.0/
-rw-r--r--	root/root	usr/share/gir-1.0/SpiceClientGLib-2.0.gir
-rw-r--r--	root/root	usr/share/gir-1.0/SpiceClientGtk-3.0.gir
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/cs/
drwxr-xr-x	root/root	usr/share/locale/cs/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/cs/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/de/
drwxr-xr-x	root/root	usr/share/locale/de/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/de/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/el/
drwxr-xr-x	root/root	usr/share/locale/el/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/el/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/es/
drwxr-xr-x	root/root	usr/share/locale/es/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/es/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/fr/
drwxr-xr-x	root/root	usr/share/locale/fr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fr/LC_MESSAGES/spice-gtk.mo
drwxr-xr-x	root/root	usr/share/locale/hu/
drwxr-xr-x	root/root	usr/share/locale/hu/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/hu/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/pl/
drwxr-xr-x	root/root	usr/share/locale/pl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pl/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/pt_BR/
drwxr-xr-x	root/root	usr/share/locale/pt_BR/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pt_BR/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/sl/
drwxr-xr-x	root/root	usr/share/locale/sl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sl/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/sr/
drwxr-xr-x	root/root	usr/share/locale/sr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sr/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/locale/sr@latin/
drwxr-xr-x	root/root	usr/share/locale/sr@latin/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sr@latin/LC_MESSAGES/phodav.mo
drwxr-xr-x	root/root	usr/share/vala/
drwxr-xr-x	root/root	usr/share/vala/vapi/
-rw-r--r--	root/root	usr/share/vala/vapi/spice-protocol.vapi
