# ts/ports/components/ntp-TS/.footprint
drwxr-xr-x	root/root	ts/
drwxr-xr-x	root/root	ts/components/
drwxr-xr-x	root/root	ts/components/ntp/
drwxr-xr-x	root/root	ts/components/ntp/etc/
drwxr-xr-x	root/root	ts/components/ntp/etc/ntp/
-rw-r--r--	root/root	ts/components/ntp/etc/ntp/ntp.conf
drwxr-xr-x	root/root	ts/components/ntp/etc/rc.d/
-rwxr-xr-x	root/root	ts/components/ntp/etc/rc.d/ntpd
-rwxr-xr-x	root/root	ts/components/ntp/etc/rc.d/ntpdate
drwxr-xr-x	root/root	ts/components/ntp/usr/
drwxr-xr-x	root/root	ts/components/ntp/usr/bin/
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntp-keygen
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntp-wait
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntpd
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntpdate
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntpdc
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntpq
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntptime
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/ntptrace
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/sntp
-rwxr-xr-x	root/root	ts/components/ntp/usr/bin/tickadj
drwxr-xr-x	root/root	ts/components/ntp/usr/lib/
