# ts/build/packages/autoruns/etc/xdg/autostart/autoruns.desktop
[Desktop Entry]
Version=1.0
Name=autoruns
Comment=Runs startup scripts for the autoruns package
TryExec=/bin/autoruns
Exec=/bin/autoruns xfce
StartupNotify=false
Terminal=false
Type=Application
Icon=ica
Categories=GTK;GNOME;X-Xfce-Toplevel;