# ts/ports/core/rsync/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/ports/
drwxr-xr-x	root/root	etc/ports/drivers/
-rwxr-xr-x	root/root	etc/ports/drivers/rsync
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/rsyncd
-rw-r--r--	root/root	etc/rsyncd.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/rsync
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/rsync.1.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/rsyncd.conf.5.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
-rw-r--r--	root/root	var/log/rsyncd.log (EMPTY)
