# ts/build/packages/autofs/build/extra/etc/systemd/system/autofs.service
[Unit]
Description=ThinStation Network AutoFS
After=profile-setup.service pkg.service
Before=display-manager.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/autofs start
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target

