# ts/ports/opt/rpcbind/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/rpcbind
drwxr-xr-x	root/root	sbin/
-rwxr-xr-x	root/root	sbin/rpcbind
-rwxr-xr-x	root/root	sbin/rpcinfo
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/rpcbind.8.gz
-rw-r--r--	root/root	usr/man/man8/rpcinfo.8.gz
