# ts/build/packages/udisks/build/extra/etc/xdg/autostart/udisks-automount.desktop
[Desktop Entry]
Version=1.0
Name=udisks-automount
Comment=Mount all available disks
Exec=/bin/udisks-automount.sh
StartupNotify=false
Terminal=false
Type=Application
Categories=GTK;GNOME;X-Xfce-Toplevel;
