# ts/build/packages/networkmanager/build/extra/lib/systemd/system/nmcreatesysconn.service
[Unit]
Description=Network Manager System Connection Creator
After=profile_setup.service
Before=NetworkManager.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/lib/thinstation/system/nmcreatesysconn
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
