# ts/build/packages/splash/etc/systemd/system/splash-start.service
[Unit]
Description=Start Splash
DefaultDependencies=no
After=pre-udev.service

[Service]
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/splash start
SyslogIdentifier=thinstation
Type=simple
KillMode=none
SendSIGKILL=no
RemainAfterExit=false
