# ts/build/packages/splash/etc/systemd/system/splash-stop.service
[Unit]
Description=Stop Splash
DefaultDependencies=no
After=systemd-udev-settle.service display-manager.service

[Service]
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/splash stop
SyslogIdentifier=thinstation
Type=oneshot
TimeoutSec=20

