# ts/build/packages/splash/etc/systemd/system/splash-start.service
[Unit]
Description=Start Splash
DefaultDependencies=no
Before=fastboot.service sysinit.target

[Service]
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/splash start
SyslogIdentifier=thinstation
Type=forking
KillMode=none
SendSIGKILL=no
