# ts/ports/contrib/acpid/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/acpi/
drwxr-xr-x	root/root	etc/acpi/events/
-rw-r--r--	root/root	etc/acpi/events/anything
-rwxr-xr-x	root/root	etc/acpi/handler.sh
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/acpi_listen
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/acpid.service
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/acpid
-rwxr-xr-x	root/root	usr/sbin/kacpimon
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/acpid/
-rw-r--r--	root/root	usr/share/doc/acpid/COPYING
-rw-r--r--	root/root	usr/share/doc/acpid/Changelog
-rw-r--r--	root/root	usr/share/doc/acpid/README
-rw-r--r--	root/root	usr/share/doc/acpid/TESTPLAN
-rw-r--r--	root/root	usr/share/doc/acpid/TODO
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/acpi_listen.8.gz
-rw-r--r--	root/root	usr/share/man/man8/acpid.8.gz
-rw-r--r--	root/root	usr/share/man/man8/kacpimon.8.gz
