# usr/ports/opt/nfs-utils/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/exports
drwxr-xr-x	root/root	etc/exports.d/
-rw-r--r--	root/root	etc/idmapd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/nfs
-rwxr-xr-x	root/root	etc/rc.d/nfsserver
drwxr-xr-x	root/root	sbin/
-r-s--x--x	root/root	sbin/mount.nfs
lrwxrwxrwx	root/root	sbin/mount.nfs4 -> mount.nfs
-rwxr-xr-x	root/root	sbin/osd_login
lrwxrwxrwx	root/root	sbin/umount.nfs -> mount.nfs
lrwxrwxrwx	root/root	sbin/umount.nfs4 -> mount.nfs
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/exports.5.gz
-rw-r--r--	root/root	usr/man/man5/nfs.5.gz
drwxr-xr-x	root/root	usr/man/man7/
-rw-r--r--	root/root	usr/man/man7/nfsd.7.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/blkmapd.8.gz
-rw-r--r--	root/root	usr/man/man8/exportfs.8.gz
-rw-r--r--	root/root	usr/man/man8/idmapd.8.gz
-rw-r--r--	root/root	usr/man/man8/mount.nfs.8.gz
-rw-r--r--	root/root	usr/man/man8/mountd.8.gz
-rw-r--r--	root/root	usr/man/man8/mountstats.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsd.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsdcltrack.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsidmap.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsiostat.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsstat.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.idmapd.8.gz -> idmapd.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.mountd.8.gz -> mountd.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.nfsd.8.gz -> nfsd.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.sm-notify.8.gz -> sm-notify.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.statd.8.gz -> statd.8.gz
-rw-r--r--	root/root	usr/man/man8/rpcdebug.8.gz
-rw-r--r--	root/root	usr/man/man8/showmount.8.gz
-rw-r--r--	root/root	usr/man/man8/sm-notify.8.gz
-rw-r--r--	root/root	usr/man/man8/statd.8.gz
-rw-r--r--	root/root	usr/man/man8/umount.nfs.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/blkmapd
-rwxr-xr-x	root/root	usr/sbin/exportfs
-rwxr-xr-x	root/root	usr/sbin/mountstats
-rwxr-xr-x	root/root	usr/sbin/nfsdcltrack
-rwxr-xr-x	root/root	usr/sbin/nfsidmap
-rwxr-xr-x	root/root	usr/sbin/nfsiostat
-rwxr-xr-x	root/root	usr/sbin/nfsstat
-rwxr-xr-x	root/root	usr/sbin/rpc.idmapd
-rwxr-xr-x	root/root	usr/sbin/rpc.mountd
-rwxr-xr-x	root/root	usr/sbin/rpc.nfsd
-rwxr-xr-x	root/root	usr/sbin/rpc.statd
-rwxr-xr-x	root/root	usr/sbin/rpcdebug
-rwxr-xr-x	root/root	usr/sbin/showmount
-rwxr-xr-x	root/root	usr/sbin/sm-notify
-rwxr-xr-x	root/root	usr/sbin/start-statd
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	nobody/root	var/lib/nfs/
-rw-r--r--	root/root	var/lib/nfs/etab (EMPTY)
-rw-r--r--	root/root	var/lib/nfs/rmtab (EMPTY)
drwxr-xr-x	root/root	var/lib/nfs/rpc_pipefs/
drwx------	nobody/root	var/lib/nfs/sm.bak/
drwx------	nobody/root	var/lib/nfs/sm/
-rw-------	nobody/root	var/lib/nfs/state (EMPTY)
drwxr-xr-x	root/root	var/lib/nfs/v4recovery/
-rw-r--r--	root/root	var/lib/nfs/xtab (EMPTY)
