# ts/build/packages/www/etc/systemd/system/www.service
[Unit]
Description=ThinStation www
After=profile-setup.service pkg.service
Before=display-manager.service lighttpd.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/www init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
