# ts/build/packages/firewall/etc/systemd/system/iptables.service
[Unit]
Description=Load Iptables Rules
After=network.target

[Service]
Type=forking
ExecStart=/sbin/iptables-restore /etc/sysconfig/iptables
TimeoutSec=0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
