# ts/5.2/packages/ntp/etc/NetworkManager/dispatcher.d/90-ntp
#!/bin/sh

. /etc/thinstation.global

INTERFACE=$1
ACTION=$2

case $ACTION in

up)
  /etc/init.d/ntpd start
;;

*)
  exit
;;

esac

