# ts/build/packages/local-install/etc/systemd/system/local-install-mount.service
[Unit]
Description=ThinStation Local Install Mount
After=profile-setup.service
Before=network-pre.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/local-install-mount
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target

