# ts/ports/components/libcanberra/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/canberra-boot
-rwxr-xr-x	root/root	usr/bin/canberra-gtk-play
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/canberra-gtk.h
-rw-r--r--	root/root	usr/include/canberra.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/gnome-settings-daemon-3.0/
drwxr-xr-x	root/root	usr/lib/gnome-settings-daemon-3.0/gtk-modules/
-rw-r--r--	root/root	usr/lib/gnome-settings-daemon-3.0/gtk-modules/canberra-gtk-module.desktop
drwxr-xr-x	root/root	usr/lib/gtk-2.0/
drwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/
-rwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/libcanberra-gtk-module.la
-rwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/libcanberra-gtk-module.so
drwxr-xr-x	root/root	usr/lib/gtk-3.0/
drwxr-xr-x	root/root	usr/lib/gtk-3.0/modules/
lrwxrwxrwx	root/root	usr/lib/gtk-3.0/modules/libcanberra-gtk-module.so -> libcanberra-gtk3-module.so
-rwxr-xr-x	root/root	usr/lib/gtk-3.0/modules/libcanberra-gtk3-module.la
-rwxr-xr-x	root/root	usr/lib/gtk-3.0/modules/libcanberra-gtk3-module.so
drwxr-xr-x	root/root	usr/lib/libcanberra-0.30/
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-alsa.la
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-alsa.so
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-gstreamer.la
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-gstreamer.so
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-multi.la
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-multi.so
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-null.la
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-null.so
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-pulse.la
-rwxr-xr-x	root/root	usr/lib/libcanberra-0.30/libcanberra-pulse.so
-rwxr-xr-x	root/root	usr/lib/libcanberra-gtk.la
lrwxrwxrwx	root/root	usr/lib/libcanberra-gtk.so -> libcanberra-gtk.so.0.1.9
lrwxrwxrwx	root/root	usr/lib/libcanberra-gtk.so.0 -> libcanberra-gtk.so.0.1.9
-rwxr-xr-x	root/root	usr/lib/libcanberra-gtk.so.0.1.9
-rwxr-xr-x	root/root	usr/lib/libcanberra-gtk3.la
lrwxrwxrwx	root/root	usr/lib/libcanberra-gtk3.so -> libcanberra-gtk3.so.0.1.9
lrwxrwxrwx	root/root	usr/lib/libcanberra-gtk3.so.0 -> libcanberra-gtk3.so.0.1.9
-rwxr-xr-x	root/root	usr/lib/libcanberra-gtk3.so.0.1.9
-rwxr-xr-x	root/root	usr/lib/libcanberra.la
lrwxrwxrwx	root/root	usr/lib/libcanberra.so -> libcanberra.so.0.2.5
lrwxrwxrwx	root/root	usr/lib/libcanberra.so.0 -> libcanberra.so.0.2.5
-rwxr-xr-x	root/root	usr/lib/libcanberra.so.0.2.5
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libcanberra-gtk.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libcanberra-gtk3.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libcanberra.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/canberra-system-bootup.service
-rw-r--r--	root/root	usr/lib/systemd/system/canberra-system-shutdown-reboot.service
-rw-r--r--	root/root	usr/lib/systemd/system/canberra-system-shutdown.service
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/gdm/
drwxr-xr-x	root/root	usr/share/gdm/autostart/
drwxr-xr-x	root/root	usr/share/gdm/autostart/LoginWindow/
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/libcanberra-ready-sound.desktop
drwxr-xr-x	root/root	usr/share/gnome/
drwxr-xr-x	root/root	usr/share/gnome/autostart/
-rw-r--r--	root/root	usr/share/gnome/autostart/libcanberra-login-sound.desktop
drwxr-xr-x	root/root	usr/share/gnome/shutdown/
-rwxr-xr-x	root/root	usr/share/gnome/shutdown/libcanberra-logout-sound.sh
drwxr-xr-x	root/root	usr/share/vala/
drwxr-xr-x	root/root	usr/share/vala/vapi/
-rw-r--r--	root/root	usr/share/vala/vapi/libcanberra-gtk.vapi
-rw-r--r--	root/root	usr/share/vala/vapi/libcanberra.vapi
