# ts/ports/opt/nfs-utils/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/exports
drwxr-xr-x	root/root	etc/exports.d/
-rw-r--r--	root/root	etc/idmapd.conf
-rw-r--r--	root/root	etc/nfsmount.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/nfs
-rwxr-xr-x	root/root	etc/rc.d/nfsclient
-rwxr-xr-x	root/root	etc/rc.d/nfsserver
-rwxr-xr-x	root/root	etc/rc.d/rpc.idmapd
-rwxr-xr-x	root/root	etc/rc.d/rpc.mountd
-rwxr-xr-x	root/root	etc/rc.d/rpc.nfsd
-rwxr-xr-x	root/root	etc/rc.d/rpc.statd
drwxr-xr-x	root/root	sbin/
-r-s--x--x	root/root	sbin/mount.nfs
lrwxrwxrwx	root/root	sbin/mount.nfs4 -> mount.nfs
-rwxr-xr-x	root/root	sbin/nfsdcltrack
lrwxrwxrwx	root/root	sbin/umount.nfs -> mount.nfs
lrwxrwxrwx	root/root	sbin/umount.nfs4 -> mount.nfs
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/nfsidmap.h
-rw-r--r--	root/root	usr/include/nfsidmap_plugin.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libnfsidmap.a
-rwxr-xr-x	root/root	usr/lib/libnfsidmap.la
lrwxrwxrwx	root/root	usr/lib/libnfsidmap.so -> libnfsidmap.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libnfsidmap.so.1 -> libnfsidmap.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libnfsidmap.so.1.0.0
drwxr-xr-x	root/root	usr/lib/libnfsidmap/
-rw-r--r--	root/root	usr/lib/libnfsidmap/nsswitch.a
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/nsswitch.la
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/nsswitch.so
-rw-r--r--	root/root	usr/lib/libnfsidmap/regex.a
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/regex.la
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/regex.so
-rw-r--r--	root/root	usr/lib/libnfsidmap/static.a
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/static.la
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/static.so
-rw-r--r--	root/root	usr/lib/libnfsidmap/umich_ldap.a
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/umich_ldap.la
-rwxr-xr-x	root/root	usr/lib/libnfsidmap/umich_ldap.so
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libnfsidmap.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system-generators/
-rwxr-xr-x	root/root	usr/lib/systemd/system-generators/nfs-server-generator
-rwxr-xr-x	root/root	usr/lib/systemd/system-generators/rpc-pipefs-generator
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/nfs-blkmap.service
-rw-r--r--	root/root	usr/lib/systemd/system/nfs-client.target
-rw-r--r--	root/root	usr/lib/systemd/system/nfs-idmapd.service
-rw-r--r--	root/root	usr/lib/systemd/system/nfs-mountd.service
-rw-r--r--	root/root	usr/lib/systemd/system/nfs-server.service
-rw-r--r--	root/root	usr/lib/systemd/system/nfs-utils.service
-rw-r--r--	root/root	usr/lib/systemd/system/nfsdcld.service
-rw-r--r--	root/root	usr/lib/systemd/system/proc-fs-nfsd.mount
-rw-r--r--	root/root	usr/lib/systemd/system/rpc-statd-notify.service
-rw-r--r--	root/root	usr/lib/systemd/system/rpc-statd.service
-rw-r--r--	root/root	usr/lib/systemd/system/rpc_pipefs.target
-rw-r--r--	root/root	usr/lib/systemd/system/var-lib-nfs-rpc_pipefs.mount
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man3/
-rw-r--r--	root/root	usr/man/man3/nfs4_uid_to_name.3.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/exports.5.gz
-rw-r--r--	root/root	usr/man/man5/idmapd.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/nfs.5.gz
-rw-r--r--	root/root	usr/man/man5/nfs.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/nfsmount.conf.5.gz
drwxr-xr-x	root/root	usr/man/man7/
-rw-r--r--	root/root	usr/man/man7/nfs.systemd.7.gz
-rw-r--r--	root/root	usr/man/man7/nfsd.7.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/blkmapd.8.gz
-rw-r--r--	root/root	usr/man/man8/exportfs.8.gz
-rw-r--r--	root/root	usr/man/man8/idmapd.8.gz
-rw-r--r--	root/root	usr/man/man8/mount.nfs.8.gz
-rw-r--r--	root/root	usr/man/man8/mountd.8.gz
-rw-r--r--	root/root	usr/man/man8/mountstats.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsconf.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsd.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsdcld.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsdclddb.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsdclnts.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsdcltrack.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsidmap.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsiostat.8.gz
-rw-r--r--	root/root	usr/man/man8/nfsstat.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.idmapd.8.gz -> idmapd.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.mountd.8.gz -> mountd.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.nfsd.8.gz -> nfsd.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.sm-notify.8.gz -> sm-notify.8.gz
lrwxrwxrwx	root/root	usr/man/man8/rpc.statd.8.gz -> statd.8.gz
-rw-r--r--	root/root	usr/man/man8/rpcdebug.8.gz
-rw-r--r--	root/root	usr/man/man8/showmount.8.gz
-rw-r--r--	root/root	usr/man/man8/sm-notify.8.gz
-rw-r--r--	root/root	usr/man/man8/statd.8.gz
-rw-r--r--	root/root	usr/man/man8/umount.nfs.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/blkmapd
-rwxr-xr-x	root/root	usr/sbin/exportfs
-rwxr-xr-x	root/root	usr/sbin/mountstats
-rwxr-xr-x	root/root	usr/sbin/nfsconf
-rwxr-xr-x	root/root	usr/sbin/nfsdcld
-rwxr-xr-x	root/root	usr/sbin/nfsdclddb
-rwxr-xr-x	root/root	usr/sbin/nfsdclnts
-rwxr-xr-x	root/root	usr/sbin/nfsidmap
-rwxr-xr-x	root/root	usr/sbin/nfsiostat
-rwxr-xr-x	root/root	usr/sbin/nfsstat
-rwxr-xr-x	root/root	usr/sbin/rpc.idmapd
-rwxr-xr-x	root/root	usr/sbin/rpc.mountd
-rwxr-xr-x	root/root	usr/sbin/rpc.nfsd
-rwxr-xr-x	root/root	usr/sbin/rpc.statd
-rwxr-xr-x	root/root	usr/sbin/rpcdebug
-rwxr-xr-x	root/root	usr/sbin/showmount
-rwxr-xr-x	root/root	usr/sbin/sm-notify
-rwxr-xr-x	root/root	usr/sbin/start-statd
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	nobody/root	var/lib/nfs/
-rw-r--r--	root/root	var/lib/nfs/etab (EMPTY)
drwxr-xr-x	root/root	var/lib/nfs/nfsdcltrack/
-rw-r--r--	root/root	var/lib/nfs/rmtab (EMPTY)
drwxr-xr-x	root/root	var/lib/nfs/rpc_pipefs/
drwx------	nobody/root	var/lib/nfs/sm.bak/
drwx------	nobody/root	var/lib/nfs/sm/
-rw-------	nobody/root	var/lib/nfs/state (EMPTY)
drwxr-xr-x	root/root	var/lib/nfs/v4recovery/
