# ts/build/packages/locale/etc/systemd/system/console-keymap.service
[Unit]
Description=Configure Locale and Keymap Settings
After=session.service pkg.service
Before=session-setup.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/console-keymap
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
