# ts/ports/components/audit/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/audisp/
-rw-r-----	root/root	etc/audisp/audisp-remote.conf
-rw-r-----	root/root	etc/audisp/audispd.conf
drwxr-x---	root/root	etc/audisp/plugins.d/
-rw-r-----	root/root	etc/audisp/plugins.d/af_unix.conf
-rw-r-----	root/root	etc/audisp/plugins.d/au-remote.conf
-rw-r-----	root/root	etc/audisp/plugins.d/audispd-zos-remote.conf
-rw-r-----	root/root	etc/audisp/plugins.d/syslog.conf
-rw-r-----	root/root	etc/audisp/zos-remote.conf
drwxr-xr-x	root/root	etc/audit/
-rw-r--r--	root/root	etc/audit/audit-stop.rules
-rw-r--r--	root/root	etc/audit/auditd.conf
-rw-r-----	root/root	etc/libaudit.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/aulast
-rwxr-xr-x	root/root	usr/bin/aulastlog
-rwxr-xr-x	root/root	usr/bin/ausyscall
-rwxr-xr-x	root/root	usr/bin/auvirt
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/auparse-defs.h
-rw-r--r--	root/root	usr/include/auparse.h
-rw-r--r--	root/root	usr/include/libaudit.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/golang/
drwxr-xr-x	root/root	usr/lib/golang/src/
drwxr-xr-x	root/root	usr/lib/golang/src/pkg/
drwxr-xr-x	root/root	usr/lib/golang/src/pkg/redhat.com/
drwxr-xr-x	root/root	usr/lib/golang/src/pkg/redhat.com/audit/
-rw-r--r--	root/root	usr/lib/golang/src/pkg/redhat.com/audit/audit.go
-rw-r--r--	root/root	usr/lib/libaudit.a
-rwxr-xr-x	root/root	usr/lib/libaudit.la
lrwxrwxrwx	root/root	usr/lib/libaudit.so -> libaudit.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libaudit.so.1 -> libaudit.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libaudit.so.1.0.0
-rw-r--r--	root/root	usr/lib/libauparse.a
-rwxr-xr-x	root/root	usr/lib/libauparse.la
lrwxrwxrwx	root/root	usr/lib/libauparse.so -> libauparse.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libauparse.so.0 -> libauparse.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libauparse.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/audit.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/auparse.pc
drwxr-xr-x	root/root	usr/lib/python3.8/
drwxr-xr-x	root/root	usr/lib/python3.8/site-packages/
drwxr-xr-x	root/root	usr/lib/python3.8/site-packages/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.8/site-packages/__pycache__/audit.cpython-38.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.8/site-packages/__pycache__/audit.cpython-38.pyc
-rwxr-xr-x	root/root	usr/lib/python3.8/site-packages/_audit.la
-rwxr-xr-x	root/root	usr/lib/python3.8/site-packages/_audit.so
-rw-r--r--	root/root	usr/lib/python3.8/site-packages/audit.py
-rw-r--r--	root/root	usr/lib/python3.8/site-packages/auparse.a
-rwxr-xr-x	root/root	usr/lib/python3.8/site-packages/auparse.la
-rwxr-xr-x	root/root	usr/lib/python3.8/site-packages/auparse.so
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/auditd.service
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/audisp-remote
-rwxr-x---	root/root	usr/sbin/audispd
-rwxr-xr-x	root/root	usr/sbin/audispd-zos-remote
-rwxr-xr-x	root/root	usr/sbin/auditctl
-rwxr-xr-x	root/root	usr/sbin/auditd
-rwxr-x---	root/root	usr/sbin/augenrules
-rwxr-xr-x	root/root	usr/sbin/aureport
-rwxr-xr-x	root/root	usr/sbin/ausearch
-rwxr-xr-x	root/root	usr/sbin/autrace
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/aclocal/
-rw-r--r--	root/root	usr/share/aclocal/audit.m4
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man3/
-rw-r--r--	root/root	usr/share/man/man3/audit_add_rule_data.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_add_watch.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_delete_rule_data.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_detect_machine.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_encode_nv_string.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_get_reply.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_get_session.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_getloginuid.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_log_acct_message.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_log_semanage_message.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_log_user_avc_message.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_log_user_comm_message.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_log_user_command.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_log_user_message.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_open.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_request_rules_list_data.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_request_signal_info.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_request_status.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_set_backlog_limit.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_set_backlog_wait_time.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_set_enabled.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_set_failure.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_set_pid.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_set_rate_limit.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_setloginuid.3.gz
-rw-r--r--	root/root	usr/share/man/man3/audit_update_watch_perms.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_add_callback.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_destroy.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_feed.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_feed_age_events.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_feed_has_data.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_find_field.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_find_field_next.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_first_field.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_first_record.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_flush_feed.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_field_int.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_field_name.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_field_num.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_field_str.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_field_type.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_filename.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_line_number.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_milli.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_node.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_num_fields.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_num_records.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_record_num.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_record_text.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_serial.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_time.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_timestamp.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_type.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_get_type_name.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_goto_field_num.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_goto_record_num.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_init.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_interpret_field.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_next_event.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_next_field.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_next_record.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_node_compare.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_normalize.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_normalize_functions.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_reset.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_set_escape_mode.3.gz
-rw-r--r--	root/root	usr/share/man/man3/auparse_timestamp_compare.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_add_expression.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_add_interpreted_item.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_add_item.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_add_regex.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_add_timestamp_item.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_add_timestamp_item_ex.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_clear.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_next_event.3.gz
-rw-r--r--	root/root	usr/share/man/man3/ausearch_set_stop.3.gz
-rw-r--r--	root/root	usr/share/man/man3/get_auditfail_action.3.gz
-rw-r--r--	root/root	usr/share/man/man3/set_aumessage_mode.3.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/audisp-remote.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/audispd.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/auditd.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/ausearch-expression.5.gz
-rw-r--r--	root/root	usr/share/man/man5/libaudit.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/zos-remote.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man7/
-rw-r--r--	root/root	usr/share/man/man7/audit.rules.7.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/audisp-remote.8.gz
-rw-r--r--	root/root	usr/share/man/man8/audispd-zos-remote.8.gz
-rw-r--r--	root/root	usr/share/man/man8/audispd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/auditctl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/auditd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/augenrules.8.gz
-rw-r--r--	root/root	usr/share/man/man8/aulast.8.gz
-rw-r--r--	root/root	usr/share/man/man8/aulastlog.8.gz
-rw-r--r--	root/root	usr/share/man/man8/aureport.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ausearch.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ausyscall.8.gz
-rw-r--r--	root/root	usr/share/man/man8/autrace.8.gz
-rw-r--r--	root/root	usr/share/man/man8/auvirt.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/audit/
