# ts/ports/opt/upower/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/UPower/
-rw-r--r--	root/root	etc/UPower/UPower.conf
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/interfaces/
-rw-r--r--	root/root	etc/dbus-1/interfaces/org.freedesktop.UPower.Device.xml
-rw-r--r--	root/root	etc/dbus-1/interfaces/org.freedesktop.UPower.KbdBacklight.xml
-rw-r--r--	root/root	etc/dbus-1/interfaces/org.freedesktop.UPower.Wakeups.xml
-rw-r--r--	root/root	etc/dbus-1/interfaces/org.freedesktop.UPower.xml
drwxr-xr-x	root/root	etc/dbus-1/system-services/
-rw-r--r--	root/root	etc/dbus-1/system-services/org.freedesktop.UPower.service
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.UPower.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/upower
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/libupower-glib/
-rw-r--r--	root/root	usr/include/libupower-glib/up-autocleanups.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-client.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-device.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-history-item.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-stats-item.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-types.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-version.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-wakeup-item.h
-rw-r--r--	root/root	usr/include/libupower-glib/up-wakeups.h
-rw-r--r--	root/root	usr/include/libupower-glib/upower.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/girepository-1.0/
-rw-r--r--	root/root	usr/lib/girepository-1.0/UPowerGlib-1.0.typelib
-rwxr-xr-x	root/root	usr/lib/libupower-glib.la
lrwxrwxrwx	root/root	usr/lib/libupower-glib.so -> libupower-glib.so.3.0.1
lrwxrwxrwx	root/root	usr/lib/libupower-glib.so.3 -> libupower-glib.so.3.0.1
-rwxr-xr-x	root/root	usr/lib/libupower-glib.so.3.0.1
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/upower-glib.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/upower.service
drwxr-xr-x	root/root	usr/lib/udev/
drwxr-xr-x	root/root	usr/lib/udev/rules.d/
-rw-r--r--	root/root	usr/lib/udev/rules.d/95-upower-csr.rules
-rw-r--r--	root/root	usr/lib/udev/rules.d/95-upower-hid.rules
-rw-r--r--	root/root	usr/lib/udev/rules.d/95-upower-hidpp.rules
-rw-r--r--	root/root	usr/lib/udev/rules.d/95-upower-wup.rules
drwxr-xr-x	root/root	usr/lib/udisks/
-rwxr-xr-x	root/root	usr/lib/udisks/upowerd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/gir-1.0/
-rw-r--r--	root/root	usr/share/gir-1.0/UPowerGlib-1.0.gir
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/fr/
drwxr-xr-x	root/root	usr/share/locale/fr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fr/LC_MESSAGES/upower.mo
drwxr-xr-x	root/root	usr/share/locale/it/
drwxr-xr-x	root/root	usr/share/locale/it/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/it/LC_MESSAGES/upower.mo
drwxr-xr-x	root/root	usr/share/locale/pl/
drwxr-xr-x	root/root	usr/share/locale/pl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pl/LC_MESSAGES/upower.mo
drwxr-xr-x	root/root	usr/share/locale/sv/
drwxr-xr-x	root/root	usr/share/locale/sv/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sv/LC_MESSAGES/upower.mo
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/upower.1.gz
drwxr-xr-x	root/root	usr/share/man/man7/
-rw-r--r--	root/root	usr/share/man/man7/UPower.7.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/upowerd.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/upower/
