# ts/5.2/packages/ica/build/install
#!/bin/bash

gzname=`basename $1`
tarname=`basename $1 .gz`
filename=`basename $1 .tar.gz`
INSTALLDIR=./packages/ica

etcs="passwd group nsswitch.conf issue"

bins="attr bash cat chgrp chmod chown cp date dd df echo ed false ln ls mkdir more mount mv \
ping ps red rm rmdir sed sh sleep tar touch true uname vi"

sbins="init init ip ss tc telinit"

ubins="awk basename cut dirname egrep expr find gawk grep head id install mktemp sort tr wc which whoami"

libs="ld-linux.so.2 libacl.so.1 libattr.so.1 libc.so.6 libdl.so.2 libhistory.so.6 libm.so.6 \
libncurses.so.5 libnss_compat.so.2 libnss_files.so.2 libproc-3.2.8.so libpthread.so.0 \
libreadline.so.6 librt.so.1 libgobject-2.0.so.0 libglib-2.0.so.0 libgmodule-2.0.so.0 \
libcap.so.2 libgthread-2.0.so.0 libnsl.so.1 libnss_nis.so.2 libncursesw.so.5"

ulibs="libpcre.so.0 libgmp.so.10 libxml2.so.2 libgstinterfaces-0.10.so.0 libffi.so.5 libz.so.1"

dirs="proc etc/udev/rules.d etc/init.d bin lib tmp dev/bus/usb usr/share/applications usr/local/lib/netscape/plugins \
usr/lib/gstreamer-0.10 usr/local/lib/firefox/plugins usr/local/lib/chrome/plugins"

clean_dirs="setupwfc lib linuxx86 nls dev bin PkgId tmp etc/nsswitch.conf etc/passwd etc/group etc/mtab etc/issue proc \
usr/local/lib/chrome usr/local/lib/firefox usr/lib/libgstreamer-0.10.so*"

for mdir in $dirs; do
	mkdir -p $INSTALLDIR/$mdir
done
for cbin in $bins ; do
	cp -Prf /bin/$cbin $INSTALLDIR/bin/.
done
for cbin in $sbins ; do
	cp -Prf /sbin/$cbin $INSTALLDIR/bin/.
done
for cbin in $ubins ; do
	cp -Prf /usr/bin/$cbin $INSTALLDIR/bin/.
done
for clib in $libs ; do
	cp /lib/$clib $INSTALLDIR/lib/.
done
for clib in $ulibs ; do
	cp /usr/lib/$clib $INSTALLDIR/lib/
done
for cetc in $etcs ; do
	cp -Prf /etc/$cetc $INSTALLDIR/etc/.
done
if [ -e /bin/uname-auto ] ; then
	cp /bin/uname-auto $INSTALLDIR/bin/.
fi

cat /proc/mounts |grep -Ev 'gvfs|mapper|media' > $INSTALLDIR/etc/mtab
mount -t proc proc $INSTALLDIR/proc

echo "exit 0" > $INSTALLDIR/bin/install_initd
echo "exit 0" > $INSTALLDIR/bin/update-rc.d
echo "exit 0" > $INSTALLDIR/bin/chkconfig
chmod +x $INSTALLDIR/bin/install_initd
chmod +x $INSTALLDIR/bin/update-rc.d
chmod +x $INSTALLDIR/bin/chkconfig

touch $INSTALLDIR/bin/firefox
touch $INSTALLDIR/bin/chrome

cp -P /usr/lib/{libgstreamer-0.10.so,libgstreamer-0.10.so.0.29.0,libgstreamer-0.10.so.0} $INSTALLDIR/usr/lib/.

 tar -C $INSTALLDIR -zxvf ./wget_tmp/$gzname
 let returnval=$?

echo -e "
Running the ICA insaller. \n\
You will need to punch these keys manually. \n\
Choose 1 to install and then accept the default location. \n\
Choose y to the next 2 questions. \n\
Choose 1 to accept the ICA license. \n\
Choose y for the next 3 questions. \n\
Choose 3 to exit. \n\
\n\
Don't worry about the errors. :) \n\
"
read -p "Press [Enter] key to start ICA installer"

chroot $INSTALLDIR ./setupwfc #< $INSTALLDIR/build/answerfile
chroot $INSTALLDIR ./etc/init.d/ctxusbd stop
umount $INSTALLDIR/proc
sleep 1

for dirs in $clean_dirs ; do
	rm -rf $INSTALLDIR/$dirs
done
mv $INSTALLDIR/usr/lib $INSTALLDIR/.
mv $INSTALLDIR/usr/share/applications $INSTALLDIR/lib/.
rm -rf $INSTALLDIR/usr
rm $INSTALLDIR/opt/Citrix/ICAClient/util/gst_play
rm $INSTALLDIR/opt/Citrix/ICAClient/util/gst_read
rm $INSTALLDIR/opt/Citrix/ICAClient/util/libgstflatstm.so
ln -sf gst_play.32 $INSTALLDIR/opt/Citrix/ICAClient/util/gst_play
ln -sf gst_read.32 $INSTALLDIR/opt/Citrix/ICAClient/util/gst_read
ln -sf libgstflatstm.32.so $INSTALLDIR/opt/Citrix/ICAClient/util/libgstflatstm.so
rm $INSTALLDIR/opt/Citrix/ICAClient/util/libgstflatstm.64.so $INSTALLDIR/opt/Citrix/ICAClient/util/gst_read.64 $INSTALLDIR/opt/Citrix/ICAClient/util/gst_play.64
find $INSTALLDIR -group sys |cut -d '/' -f4- > $INSTALLDIR/etc/ICA-Perms
mkdir -p $INSTALLDIR/lib/pcsc/drivers
cp -Prf $INSTALLDIR/build/extra/* $INSTALLDIR/.
