# ts/ports/components/tigervnc/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/tigervnc
drwxr-xr-x	root/root	etc/tigervnc/
-rw-r--r--	root/root	etc/tigervnc/vncserver-config-defaults
-rw-r--r--	root/root	etc/tigervnc/vncserver-config-mandatory
-rw-r--r--	root/root	etc/tigervnc/vncserver.users
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/Xvnc
-rwxr-xr-x	root/root	usr/bin/vncconfig
-rwxr-xr-x	root/root	usr/bin/vncpasswd
-rwxr-xr-x	root/root	usr/bin/vncviewer
-rwxr-xr-x	root/root	usr/bin/x0vncserver
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/vncserver@.service
drwxr-xr-x	root/root	usr/lib/xorg/
drwxr-xr-x	root/root	usr/lib/xorg/modules/
drwxr-xr-x	root/root	usr/lib/xorg/modules/extensions/
-rwxr-xr-x	root/root	usr/lib/xorg/modules/extensions/libvnc.la
-rwxr-xr-x	root/root	usr/lib/xorg/modules/extensions/libvnc.so
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/vncserver
-rwxr-xr-x	root/root	usr/libexec/vncsession-start
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/vncsession
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/applications/
-rw-r--r--	root/root	usr/share/applications/vncviewer.desktop
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/tigervnc/
-rw-r--r--	root/root	usr/share/doc/tigervnc/HOWTO.md
-rw-r--r--	root/root	usr/share/doc/tigervnc/LICENCE.TXT
-rw-r--r--	root/root	usr/share/doc/tigervnc/README.rst
drwxr-xr-x	root/root	usr/share/icons/
drwxr-xr-x	root/root	usr/share/icons/hicolor/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/16x16/apps/tigervnc.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/22x22/apps/tigervnc.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/24x24/apps/tigervnc.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/32x32/apps/tigervnc.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/48x48/
drwxr-xr-x	root/root	usr/share/icons/hicolor/48x48/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/48x48/apps/tigervnc.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/scalable/apps/tigervnc.svg
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/bg/
drwxr-xr-x	root/root	usr/share/locale/bg/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/bg/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/cs/
drwxr-xr-x	root/root	usr/share/locale/cs/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/cs/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/da/
drwxr-xr-x	root/root	usr/share/locale/da/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/da/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/de/
drwxr-xr-x	root/root	usr/share/locale/de/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/de/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/el/
drwxr-xr-x	root/root	usr/share/locale/el/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/el/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/eo/
drwxr-xr-x	root/root	usr/share/locale/eo/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/eo/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/es/
drwxr-xr-x	root/root	usr/share/locale/es/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/es/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/fi/
drwxr-xr-x	root/root	usr/share/locale/fi/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fi/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/fr/
drwxr-xr-x	root/root	usr/share/locale/fr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fr/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/fur/
drwxr-xr-x	root/root	usr/share/locale/fur/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fur/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/he/
drwxr-xr-x	root/root	usr/share/locale/he/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/he/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/hu/
drwxr-xr-x	root/root	usr/share/locale/hu/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/hu/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/id/
drwxr-xr-x	root/root	usr/share/locale/id/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/id/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/it/
drwxr-xr-x	root/root	usr/share/locale/it/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/it/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/ko/
drwxr-xr-x	root/root	usr/share/locale/ko/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ko/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/nl/
drwxr-xr-x	root/root	usr/share/locale/nl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/nl/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/pl/
drwxr-xr-x	root/root	usr/share/locale/pl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pl/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/pt_BR/
drwxr-xr-x	root/root	usr/share/locale/pt_BR/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pt_BR/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/ru/
drwxr-xr-x	root/root	usr/share/locale/ru/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ru/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/sk/
drwxr-xr-x	root/root	usr/share/locale/sk/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sk/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/sr/
drwxr-xr-x	root/root	usr/share/locale/sr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sr/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/sv/
drwxr-xr-x	root/root	usr/share/locale/sv/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sv/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/tr/
drwxr-xr-x	root/root	usr/share/locale/tr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/tr/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/uk/
drwxr-xr-x	root/root	usr/share/locale/uk/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/uk/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/vi/
drwxr-xr-x	root/root	usr/share/locale/vi/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/vi/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/zh_CN/
drwxr-xr-x	root/root	usr/share/locale/zh_CN/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/zh_CN/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/locale/zh_TW/
drwxr-xr-x	root/root	usr/share/locale/zh_TW/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/zh_TW/LC_MESSAGES/tigervnc.mo
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/Xvnc.1.gz
-rw-r--r--	root/root	usr/share/man/man1/vncconfig.1.gz
-rw-r--r--	root/root	usr/share/man/man1/vncpasswd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/vncviewer.1.gz
-rw-r--r--	root/root	usr/share/man/man1/x0vncserver.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/vncserver.8.gz
-rw-r--r--	root/root	usr/share/man/man8/vncsession.8.gz
