# ts/ports/components/open-vm-tools/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/pam.d/
-rwxr-xr-x	root/root	etc/pam.d/vmtoolsd
-rw-r--r--	root/root	etc/pam.d/vmware-guestd
drwxr-xr-x	root/root	etc/vmware-tools/
-rwxr-xr-x	root/root	etc/vmware-tools/poweroff-vm-default
-rwxr-xr-x	root/root	etc/vmware-tools/poweron-vm-default
-rwxr-xr-x	root/root	etc/vmware-tools/resume-vm-default
drwxr-xr-x	root/root	etc/vmware-tools/scripts/
drwxr-xr-x	root/root	etc/vmware-tools/scripts/vmware/
-rwxr-xr-x	root/root	etc/vmware-tools/scripts/vmware/network
-rwxr-xr-x	root/root	etc/vmware-tools/statechange.subr
-rwxr-xr-x	root/root	etc/vmware-tools/suspend-vm-default
-rw-r--r--	root/root	etc/vmware-tools/tools.conf
-rw-r--r--	root/root	etc/vmware-tools/tools.conf.example
-rw-r--r--	root/root	etc/vmware-tools/xautostart.conf
drwxr-xr-x	root/root	etc/xdg/
drwxr-xr-x	root/root	etc/xdg/autostart/
-rw-r--r--	root/root	etc/xdg/autostart/vmware-user.desktop
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/vm-support
-rwxr-xr-x	root/root	usr/bin/vmhgfs-fuse
-rwxr-xr-x	root/root	usr/bin/vmtoolsd
-rwxr-xr-x	root/root	usr/bin/vmware-checkvm
-rwxr-xr-x	root/root	usr/bin/vmware-hgfsclient
-rwxr-xr-x	root/root	usr/bin/vmware-namespace-cmd
-rwxr-xr-x	root/root	usr/bin/vmware-rpctool
-rwxr-xr-x	root/root	usr/bin/vmware-toolbox-cmd
lrwxrwxrwx	root/root	usr/bin/vmware-user -> vmware-user-suid-wrapper
-rwxr-xr-x	root/root	usr/bin/vmware-user-suid-wrapper
-rwxr-xr-x	root/root	usr/bin/vmware-vmblock-fuse
-rwxr-xr-x	root/root	usr/bin/vmware-xferlogs
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/libDeployPkg/
-rw-r--r--	root/root	usr/include/libDeployPkg/deploypkg.h
-rw-r--r--	root/root	usr/include/libDeployPkg/guestcust-events.h
-rw-r--r--	root/root	usr/include/libDeployPkg/guestrpc.h
-rw-r--r--	root/root	usr/include/libDeployPkg/imgcust-api.h
-rw-r--r--	root/root	usr/include/libDeployPkg/includeCheck.h
-rw-r--r--	root/root	usr/include/libDeployPkg/linuxDeployment.h
-rw-r--r--	root/root	usr/include/libDeployPkg/log.h
-rw-r--r--	root/root	usr/include/libDeployPkg/process.h
-rw-r--r--	root/root	usr/include/libDeployPkg/vm_basic_types.h
drwxr-xr-x	root/root	usr/include/vmGuestLib/
-rw-r--r--	root/root	usr/include/vmGuestLib/includeCheck.h
-rw-r--r--	root/root	usr/include/vmGuestLib/vmGuestLib.h
-rw-r--r--	root/root	usr/include/vmGuestLib/vmSessionId.h
-rw-r--r--	root/root	usr/include/vmGuestLib/vm_basic_types.h
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/libDeployPkg.la
lrwxrwxrwx	root/root	usr/lib/libDeployPkg.so -> libDeployPkg.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libDeployPkg.so.0 -> libDeployPkg.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libDeployPkg.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libguestlib.la
lrwxrwxrwx	root/root	usr/lib/libguestlib.so -> libguestlib.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libguestlib.so.0 -> libguestlib.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libguestlib.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libhgfs.la
lrwxrwxrwx	root/root	usr/lib/libhgfs.so -> libhgfs.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libhgfs.so.0 -> libhgfs.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libhgfs.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libvmtools.la
lrwxrwxrwx	root/root	usr/lib/libvmtools.so -> libvmtools.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libvmtools.so.0 -> libvmtools.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libvmtools.so.0.0.0
drwxr-xr-x	root/root	usr/lib/open-vm-tools/
drwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/
drwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/common/
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/common/libhgfsServer.la
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/common/libhgfsServer.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/common/libvix.la
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/common/libvix.so
drwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/libappInfo.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/libdeployPkgPlugin.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/libresolutionKMS.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
drwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmusr/
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmusr/libdesktopEvents.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmusr/libdndcp.so
-rwxr-xr-x	root/root	usr/lib/open-vm-tools/plugins/vmusr/libresolutionSet.so
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libDeployPkg.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/vmguestlib.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/vmtoolsd.service
-rw-r--r--	root/root	usr/lib/systemd/system/vmware-vmblock-fuse.service
drwxr-xr-x	root/root	usr/lib/udev/
drwxr-xr-x	root/root	usr/lib/udev/rules.d/
-rw-r--r--	root/root	usr/lib/udev/rules.d/99-vmware-scsi-udev.rules
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/mount.vmhgfs
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/open-vm-tools/
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/de/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/de/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/de/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/en/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/en/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/en/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/es/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/es/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/es/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/fr/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/fr/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/fr/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/it/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/it/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/it/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/ja/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/ja/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/ja/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/ko/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/ko/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/ko/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/zh_CN/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/zh_CN/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/zh_CN/vmtoolsd.vmsg
drwxr-xr-x	root/root	usr/share/open-vm-tools/messages/zh_TW/
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/zh_TW/toolboxcmd.vmsg
-rw-r--r--	root/root	usr/share/open-vm-tools/messages/zh_TW/vmtoolsd.vmsg
