# ts/build/packages/octoprint/build/extra/etc/systemd/system/octoprint.service
[Unit]
    Description=OctoPrint service
    After=pkg.service

[Service]
    User=octoprint
    ExecStart=/bin/python -m octoprint

[Install]
    WantedBy=multi-user.target
