# ts/ports/components/lighttpd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/lighttpd/
-rw-r--r--	root/root	etc/lighttpd/lighttpd.conf
drwxr-xr-x	root/root	etc/logrotate.d/
-rw-r--r--	root/root	etc/logrotate.d/lighttpd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/lighttpd/
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_access.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_access.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_accesslog.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_accesslog.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_alias.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_alias.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_auth.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_auth.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_file.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_file.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_gssapi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_gssapi.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_ldap.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_ldap.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_mysql.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_authn_mysql.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_cgi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_cgi.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_cml.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_cml.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_compress.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_compress.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_deflate.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_deflate.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_dirlisting.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_dirlisting.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_evasive.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_evasive.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_evhost.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_evhost.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_expire.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_expire.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_extforward.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_extforward.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_fastcgi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_fastcgi.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_flv_streaming.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_flv_streaming.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_indexfile.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_indexfile.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_magnet.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_magnet.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_mysql_vhost.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_mysql_vhost.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_openssl.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_openssl.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_proxy.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_proxy.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_redirect.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_redirect.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rewrite.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rewrite.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rrdtool.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_rrdtool.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_scgi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_scgi.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_secdownload.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_secdownload.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_setenv.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_setenv.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_simple_vhost.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_simple_vhost.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_sockproxy.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_sockproxy.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_ssi.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_ssi.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_staticfile.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_staticfile.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_status.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_status.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_trigger_b4_dl.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_trigger_b4_dl.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_uploadprogress.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_uploadprogress.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_userdir.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_userdir.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_usertrack.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_usertrack.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb_ldap.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb_ldap.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb_mysql.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_vhostdb_mysql.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_webdav.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_webdav.so
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_wstunnel.la
-rwxr-xr-x	root/root	usr/lib/lighttpd/mod_wstunnel.so
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/lighttpd.service
drwxr-xr-x	root/root	usr/lib/tmpfiles.d/
-rw-r--r--	root/root	usr/lib/tmpfiles.d/lighttpd.conf
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/lighttpd
-rwxr-xr-x	root/root	usr/sbin/lighttpd-angel
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/lighttpd/
drwxr-xr-x	root/root	usr/share/doc/lighttpd/config/
drwxr-xr-x	root/root	usr/share/doc/lighttpd/config/conf.d/
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/access_log.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/auth.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/cgi.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/cml.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/compress.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/debug.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/dirlisting.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/evhost.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/expire.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/fastcgi.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/geoip.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/magnet.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/mime.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/mod.template
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/mysql_vhost.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/proxy.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/rrdtool.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/scgi.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/secdownload.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/simple_vhost.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/ssi.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/status.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/trigger_b4_dl.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/userdir.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/conf.d/webdav.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/lighttpd.conf
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/modules.conf
drwxr-xr-x	root/root	usr/share/doc/lighttpd/config/vhosts.d/
-rw-r--r--	root/root	usr/share/doc/lighttpd/config/vhosts.d/vhosts.template
drwxr-xr-x	root/root	usr/share/licenses/
drwxr-xr-x	root/root	usr/share/licenses/lighttpd/
-rw-r--r--	root/root	usr/share/licenses/lighttpd/COPYING
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/lighttpd-angel.8.gz
-rw-r--r--	root/root	usr/share/man/man8/lighttpd.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/www/
drwx------	lighttpd/lighttpd	var/www/lighttpd/
-rw-r--r--	lighttpd/lighttpd	var/www/lighttpd/favicon.ico
-rw-r--r--	lighttpd/lighttpd	var/www/lighttpd/index.html
-rw-r--r--	lighttpd/lighttpd	var/www/lighttpd/light_button.png
-rw-r--r--	lighttpd/lighttpd	var/www/lighttpd/light_logo.png
-rw-r--r--	lighttpd/lighttpd	var/www/lighttpd/poweredby.png
