# ts/build/packages/base/etc/systemd/system/crontab.service
[Unit]
Description=ThinStation Crontab Initialize
After=profile-setup.service net.profile-setup.service
Before=crond.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/crontab
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
