# ts/build/packages/base/etc/systemd/system/fastboot.service
[Unit]
Description=Mount Fastboot Data
DefaultDependencies=no
After=systemd-udev-settle.service
Before=sysinit.target dbus.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/fastboot
SyslogIdentifier=thinstation
