# ts/build/packages/base/etc/systemd/system/mounts.service
[Unit]
Description=ThinStation Local Mounts
DefaultDependencies=no
Conflicts=shutdown.target
After=debug.service
After=systemd-udev-settle.service
Before=local-fs-pre.target local-fs.target shutdown.target
Wants=local-fs-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/mounts
SyslogIdentifier=thinstation
