# ts/build/packages/base/etc/systemd/system/post-udev.service
[Unit]
Description=ThinStation Post Udev Quirks
After=systemd-udev-settle.service
Before=display-manager.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/post-udev
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
