# ts/build/packages/base/etc/systemd/system/rtc_cmos.service
[Unit]
Description=Configure System Clock Properties
DefaultDependencies=no
After=systemd-udev-settle.service
Before=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/rtc_cmos
SyslogIdentifier=thinstation
