# ts/build/packages/base/etc/systemd/system/vmtouch.service
[Unit]
Description=Portable file system cache diagnostics and control

[Service]
Type=oneshot
ExecStart=/bin/vmtouch -ft /lib
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
