# usr/ports/core/rc/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/inittab
-rwxr-xr-x	root/root	etc/rc
-rw-r--r--	root/root	etc/rc.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.fix
-rwxr-xr-x	root/root	etc/rc.local
-rwxr-xr-x	root/root	etc/rc.modules
-rwxr-xr-x	root/root	etc/rc.multi
-rwxr-xr-x	root/root	etc/rc.shutdown
-rwxr-xr-x	root/root	etc/rc.single
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/urandom/
-rw-r-----	root/root	var/lib/urandom/seed (EMPTY)
drwxr-xr-x	root/root	var/log/
-rw-r-----	root/root	var/log/boot (EMPTY)
