# ts/ports/core/openssh/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/sshd
drwxr-xr-x	root/root	etc/rc.d/
-rw-r--r--	root/root	etc/rc.d/sshd
drwxr-xr-x	root/root	etc/ssh/
-rw-r--r--	root/root	etc/ssh/moduli
-rw-r--r--	root/root	etc/ssh/ssh_config
-rw-r--r--	root/root	etc/ssh/sshd_config
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/findssl.sh
-rwxr-xr-x	root/root	usr/bin/scp
-rwxr-xr-x	root/root	usr/bin/sftp
-rwxr-xr-x	root/root	usr/bin/ssh
-rwxr-xr-x	root/root	usr/bin/ssh-add
-rwxr-xr-x	root/root	usr/bin/ssh-agent
-rwxr-xr-x	root/root	usr/bin/ssh-copy-id
-rwxr-xr-x	root/root	usr/bin/ssh-keygen
-rwxr-xr-x	root/root	usr/bin/ssh-keyscan
-rwxr-xr-x	root/root	usr/bin/sshd
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/ssh/
-rwxr-xr-x	root/root	usr/lib/ssh/sftp-server
-rws--x--x	root/root	usr/lib/ssh/ssh-keysign
-rwxr-xr-x	root/root	usr/lib/ssh/ssh-pkcs11-helper
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/sshd.service
-rw-r--r--	root/root	usr/lib/systemd/system/sshd.socket
-rw-r--r--	root/root	usr/lib/systemd/system/sshd@.service
-rw-r--r--	root/root	usr/lib/systemd/system/sshdgenkeys.service
drwxr-xr-x	root/root	usr/lib/tmpfiles.d/
-rw-r--r--	root/root	usr/lib/tmpfiles.d/sshd.conf
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/licenses/
drwxr-xr-x	root/root	usr/share/licenses/openssh/
-rw-r--r--	root/root	usr/share/licenses/openssh/LICENCE
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/scp.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sftp.1.gz
lrwxrwxrwx	root/root	usr/share/man/man1/slogin.1.gz -> ssh.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ssh-add.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ssh-agent.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ssh-copy-id.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ssh-keygen.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ssh-keyscan.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ssh.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/moduli.5.gz
-rw-r--r--	root/root	usr/share/man/man5/ssh_config.5.gz
-rw-r--r--	root/root	usr/share/man/man5/sshd_config.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/sftp-server.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ssh-keysign.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ssh-pkcs11-helper.8.gz
-rw-r--r--	root/root	usr/share/man/man8/sshd.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/empty/
