# ts/build/packages/samba-base/etc/systemd/system/nmb.service.d/overrides.conf
[Service]
EnvironmentFile=/etc/thinstation.env
ExecStartPre=/etc/init.d/samba-base init
