# ts/build/packages/spice-vdagent/build/extra/etc/systemd/system/9p-qemu-mounts.service
[Unit]
Description=Mount shared folder between qemu host and guest
ConditionVirtualization=qemu
After=post-udev.service

[Service]
Type=oneshot
ExecStart=/bin/9p-mounts
RemainAfterExit=yes


[Install]
WantedBy=local-fs.target.wants
