# ts/build/packages/spice-vdagent/build/extra/etc/systemd/system/qemu-mounts.service
[Unit]
Description=Mount shared folder between qemu host and guest
After=post-udev.service

[Service]
Type=oneshot
ExecStart=/bin/qemu-mounts
RemainAfterExit=yes


[Install]
WantedBy=local-fs.target.wants
