# ts/build/packages/hdd-spindown/etc/systemd/system/hdd-spindown_force-init.service
[Unit]
Description=ThinStation hdd-spindown_force-init
After=profile-setup.service pkg.service
Before=display-manager.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/hdd-spindown_force-init init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
