# ts/ports/components/lightdm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.DisplayManager.conf
drwxr-xr-x	root/root	etc/lightdm/
-rwxr-xr-x	root/root	etc/lightdm/Xsession
-rw-r--r--	root/root	etc/lightdm/keys.conf
-rw-r--r--	root/root	etc/lightdm/lightdm.conf
-rw-r--r--	root/root	etc/lightdm/users.conf
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/lightdm
-rw-r--r--	root/root	etc/pam.d/lightdm-autologin
-rw-r--r--	root/root	etc/pam.d/lightdm-greeter
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/lightdm
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/dm-tool
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/lightdm-gobject-1/
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm.h
drwxr-xr-x	root/root	usr/include/lightdm-gobject-1/lightdm/
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm/greeter.h
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm/language.h
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm/layout.h
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm/power.h
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm/session.h
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm/system.h
-rw-r--r--	root/root	usr/include/lightdm-gobject-1/lightdm/user.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/girepository-1.0/
-rw-r--r--	root/root	usr/lib/girepository-1.0/LightDM-1.typelib
-rwxr-xr-x	root/root	usr/lib/liblightdm-gobject-1.la
lrwxrwxrwx	root/root	usr/lib/liblightdm-gobject-1.so -> liblightdm-gobject-1.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/liblightdm-gobject-1.so.0 -> liblightdm-gobject-1.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/liblightdm-gobject-1.so.0.0.0
drwxr-xr-x	root/root	usr/lib/lightdm/
-rwxr-xr-x	root/root	usr/lib/lightdm/lightdm-guest-session
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/liblightdm-gobject-1.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/lightdm
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/accountsservice/
drwxr-xr-x	root/root	usr/share/accountsservice/interfaces/
lrwxrwxrwx	root/root	usr/share/accountsservice/interfaces/org.freedesktop.DisplayManager.AccountsService.xml -> ../../dbus-1/interfaces/org.freedesktop.DisplayManager.AccountsService.xml
drwxr-xr-x	root/root	usr/share/bash-completion/
drwxr-xr-x	root/root	usr/share/bash-completion/completions/
-rw-r--r--	root/root	usr/share/bash-completion/completions/dm-tool
-rw-r--r--	root/root	usr/share/bash-completion/completions/lightdm
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/interfaces/
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.DisplayManager.AccountsService.xml
drwxr-xr-x	root/root	usr/share/gir-1.0/
-rw-r--r--	root/root	usr/share/gir-1.0/LightDM-1.gir
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/dm-tool.1.gz
-rw-r--r--	root/root	usr/share/man/man1/lightdm.1.gz
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/org.freedesktop.DisplayManager.AccountsService.policy
drwxr-x---	root/root	usr/share/polkit-1/rules.d/
-rw-r--r--	root/root	usr/share/polkit-1/rules.d/lightdm.rules
drwxr-xr-x	root/root	usr/share/vala/
drwxr-xr-x	root/root	usr/share/vala/vapi/
-rw-r--r--	root/root	usr/share/vala/vapi/liblightdm-gobject-1.deps
-rw-r--r--	root/root	usr/share/vala/vapi/liblightdm-gobject-1.vapi
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/cache/
drwxr-xr-t	root/root	var/cache/lightdm/
drwxr-xr-x	root/root	var/lib/
drwxrwx--T	lightdm/lightdm	var/lib/lightdm-data/
drwxrwx--T	lightdm/lightdm	var/lib/lightdm/
drwxr-xr-x	root/root	var/log/
drwx--x--x	root/lightdm	var/log/lightdm/
