# ts/ports/gnome/networkmanager/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/NetworkManager/
drwxr-xr-x	root/root	etc/NetworkManager/dispatcher.d/
drwxr-xr-x	root/root	etc/NetworkManager/system-connections/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/nm-avahi-autoipd.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/nm-dhcp-client.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/nm-dispatcher.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.NetworkManager.conf
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/77-nm-olpc-mesh.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/nm-online
-rwxr-xr-x	root/root	usr/bin/nm-tool
-rwxr-xr-x	root/root	usr/bin/nmcli
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/NetworkManager/
-rw-r--r--	root/root	usr/include/NetworkManager/NetworkManager.h
-rw-r--r--	root/root	usr/include/NetworkManager/NetworkManagerVPN.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-connection.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-8021x.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-adsl.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bluetooth.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bond.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bridge-port.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bridge.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-cdma.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-connection.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-gsm.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-infiniband.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ip4-config.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ip6-config.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-olpc-mesh.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ppp.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-pppoe.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-serial.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-vlan.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-vpn.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wimax.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wired.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wireless-security.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wireless.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-utils-enum-types.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-utils.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-version.h
drwxr-xr-x	root/root	usr/include/libnm-glib/
-rw-r--r--	root/root	usr/include/libnm-glib/libnm_glib.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-access-point.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-active-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-client.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-adsl.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-bond.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-bridge.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-bt.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-ethernet.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-infiniband.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-modem.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-olpc-mesh.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-vlan.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-wifi.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-wimax.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-dhcp4-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-dhcp6-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-glib-enum-types.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-ip4-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-ip6-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-object.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-remote-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-remote-settings.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-secret-agent.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-types.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-enum-types.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin-ui-interface.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin-utils.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-wimax-nsp.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/NetworkManager/
-rwxr-xr-x	root/root	usr/lib/libnm-glib-vpn.la
lrwxrwxrwx	root/root	usr/lib/libnm-glib-vpn.so -> libnm-glib-vpn.so.1.1.0
lrwxrwxrwx	root/root	usr/lib/libnm-glib-vpn.so.1 -> libnm-glib-vpn.so.1.1.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib-vpn.so.1.1.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib.la
lrwxrwxrwx	root/root	usr/lib/libnm-glib.so -> libnm-glib.so.4.5.0
lrwxrwxrwx	root/root	usr/lib/libnm-glib.so.4 -> libnm-glib.so.4.5.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib.so.4.5.0
-rwxr-xr-x	root/root	usr/lib/libnm-util.la
lrwxrwxrwx	root/root	usr/lib/libnm-util.so -> libnm-util.so.2.4.0
lrwxrwxrwx	root/root	usr/lib/libnm-util.so.2 -> libnm-util.so.2.4.0
-rwxr-xr-x	root/root	usr/lib/libnm-util.so.2.4.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/NetworkManager.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-glib-vpn.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-glib.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-util.pc
drwxr-xr-x	root/root	usr/lib/pppd/
drwxr-xr-x	root/root	usr/lib/pppd/2.4.5/
-rwxr-xr-x	root/root	usr/lib/pppd/2.4.5/nm-pppd-plugin.la
-rwxr-xr-x	root/root	usr/lib/pppd/2.4.5/nm-pppd-plugin.so
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/NetworkManager-dispatcher.service
-rw-r--r--	root/root	usr/lib/systemd/system/NetworkManager-wait-online.service
-rw-r--r--	root/root	usr/lib/systemd/system/NetworkManager.service
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/nm-avahi-autoipd.action
-rwxr-xr-x	root/root	usr/libexec/nm-dhcp-client.action
-rwxr-xr-x	root/root	usr/libexec/nm-dispatcher.action
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/NetworkManager
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/bash-completion/
drwxr-xr-x	root/root	usr/share/bash-completion/completions/
-rw-r--r--	root/root	usr/share/bash-completion/completions/nmcli
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.NetworkManager.service
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
drwxr-xr-x	root/root	usr/share/gtk-doc/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/NetworkManager/
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/NetworkManager.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/ix01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/license.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/ref-dbus.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/ref-migrating.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/ref-settings.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/secrets-flags.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/spec.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/NetworkManager/up.png
drwxr-xr-x	root/root	usr/share/gtk-doc/html/libnm-glib/
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMAccessPoint.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMActiveConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMClient.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDHCP4Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDHCP6Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDevice.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceAdsl.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceBond.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceBridge.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceBt.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceEthernet.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceInfiniband.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceModem.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceOlpcMesh.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceVlan.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceWifi.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceWimax.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMIP4Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMIP6Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMObject.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMRemoteConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMRemoteSettings.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMSecretAgent.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMVPNConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMWimaxNsp.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/annotation-glossary.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/api-index-full.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/ch02.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/index.sgml
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMVPNPlugin.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMVpnPluginUiWidget.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-nm-types.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-nm-vpn-plugin-utils.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/object-tree.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/ref-overview.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/up.png
drwxr-xr-x	root/root	usr/share/gtk-doc/html/libnm-util/
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSetting.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSetting8021x.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingAdsl.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingBluetooth.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingBond.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingBridge.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingBridgePort.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingCdma.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingGsm.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingIP4Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingIP6Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingInfiniband.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingOlpcMesh.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingPPP.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingPPPOE.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingSerial.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingVPN.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingVlan.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingWimax.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingWired.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingWireless.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingWirelessSecurity.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/annotation-glossary.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/api-index-full.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/ch01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/index.sgml
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/libnm-util-nm-utils.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/libnm-util.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/object-tree.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/up.png
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/nm-online.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nm-tool.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nmcli.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/NetworkManager.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/nm-system-settings.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/NetworkManager.8.gz
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwx------	root/root	var/lib/NetworkManager/
drwxr-xr-x	root/root	var/run/
drwx------	root/root	var/run/NetworkManager/
