# ts/ports/gnome/networkmanager/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/NetworkManager/
-rw-r--r--	root/root	etc/NetworkManager/NetworkManager.conf
drwxr-xr-x	root/root	etc/NetworkManager/dispatcher.d/
drwxr-xr-x	root/root	etc/NetworkManager/system-connections/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/nm-avahi-autoipd.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/nm-dhcp-client.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/nm-dispatcher.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.NetworkManager.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/networkmanager
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/77-nm-olpc-mesh.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/nm-online
-rwxr-xr-x	root/root	usr/bin/nm-tool
-rwxr-xr-x	root/root	usr/bin/nmcli
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/NetworkManager/
-rw-r--r--	root/root	usr/include/NetworkManager/NetworkManager.h
-rw-r--r--	root/root	usr/include/NetworkManager/NetworkManagerVPN.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-connection.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-8021x.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-adsl.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bluetooth.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bond.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bridge-port.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bridge.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-cdma.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-connection.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-gsm.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-infiniband.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ip4-config.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ip6-config.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-olpc-mesh.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ppp.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-pppoe.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-serial.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-vlan.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-vpn.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wimax.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wired.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wireless-security.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wireless.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-utils-enum-types.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-utils.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-version.h
drwxr-xr-x	root/root	usr/include/libnm-glib/
-rw-r--r--	root/root	usr/include/libnm-glib/libnm_glib.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-access-point.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-active-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-client.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-adsl.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-bond.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-bridge.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-bt.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-ethernet.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-infiniband.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-modem.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-olpc-mesh.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-vlan.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-wifi.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-wimax.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-dhcp4-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-dhcp6-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-glib-enum-types.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-ip4-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-ip6-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-object.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-remote-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-remote-settings.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-secret-agent.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-types.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-enum-types.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin-ui-interface.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin-utils.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-wimax-nsp.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/NetworkManager/
-rwxr-xr-x	root/root	usr/lib/libnm-glib-vpn.la
lrwxrwxrwx	root/root	usr/lib/libnm-glib-vpn.so -> libnm-glib-vpn.so.1.1.0
lrwxrwxrwx	root/root	usr/lib/libnm-glib-vpn.so.1 -> libnm-glib-vpn.so.1.1.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib-vpn.so.1.1.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib.la
lrwxrwxrwx	root/root	usr/lib/libnm-glib.so -> libnm-glib.so.4.6.0
lrwxrwxrwx	root/root	usr/lib/libnm-glib.so.4 -> libnm-glib.so.4.6.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib.so.4.6.0
-rwxr-xr-x	root/root	usr/lib/libnm-util.la
lrwxrwxrwx	root/root	usr/lib/libnm-util.so -> libnm-util.so.2.4.0
lrwxrwxrwx	root/root	usr/lib/libnm-util.so.2 -> libnm-util.so.2.4.0
-rwxr-xr-x	root/root	usr/lib/libnm-util.so.2.4.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/NetworkManager.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-glib-vpn.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-glib.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-util.pc
drwxr-xr-x	root/root	usr/lib/pppd/
drwxr-xr-x	root/root	usr/lib/pppd/2.4.7/
-rwxr-xr-x	root/root	usr/lib/pppd/2.4.7/nm-pppd-plugin.la
-rwxr-xr-x	root/root	usr/lib/pppd/2.4.7/nm-pppd-plugin.so
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/nm-avahi-autoipd.action
-rwxr-xr-x	root/root	usr/libexec/nm-dhcp-client.action
-rwxr-xr-x	root/root	usr/libexec/nm-dispatcher.action
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/nm-online.1.gz
-rw-r--r--	root/root	usr/man/man1/nm-tool.1.gz
-rw-r--r--	root/root	usr/man/man1/nmcli.1.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/NetworkManager.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/nm-system-settings.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/NetworkManager.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/NetworkManager
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwx------	root/root	var/lib/NetworkManager/
drwxr-xr-x	root/root	var/run/
drwx------	root/root	var/run/NetworkManager/
