# ts/ports/components/systemd/81-crux.rules
#
# 81-crux.rules
# 

KERNEL=="fuse", ACTION=="add", RUN+="/bin/mount -t fusectl fusectl /sys/fs/fuse/connections"

