# ts/build/packages/base/etc/inittab
# /etc/inittab init(8) configuration for BusyBox

tty1::sysinit:-/etc/thinstation.init

tty1::respawn:-/sbin/agetty -c --noclear --nohints --noissue --nohostname -a tsuser tty1 115200 linux
tty2::askfirst:/sbin/getty tty2 115200 linux
tty3::askfirst:/sbin/getty tty3 115200 linux

tty1::shutdown:/bin/shutdown

::ctrlaltdel:/sbin/reboot
