# ts/build/packages/template/etc/systemd/system/your_start_up_script.service
[Unit]
Description=ThinStation your_start_up_script
After=profile-setup.service pkg.service
Before=display-manager.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/your_start_up_script init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
