# ts/ports/components/ndctl/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/ndctl/
drwxr-xr-x	root/root	etc/ndctl/keys/
-rw-r--r--	root/root	etc/ndctl/keys/keys.readme
-rw-r--r--	root/root	etc/ndctl/monitor.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/daxctl
-rwxr-xr-x	root/root	usr/bin/ndctl
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/daxctl/
-rw-r--r--	root/root	usr/include/daxctl/libdaxctl.h
drwxr-xr-x	root/root	usr/include/ndctl/
-rw-r--r--	root/root	usr/include/ndctl/libdaxctl.h
-rw-r--r--	root/root	usr/include/ndctl/libndctl.h
-rw-r--r--	root/root	usr/include/ndctl/ndctl.h
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/libdaxctl.la
lrwxrwxrwx	root/root	usr/lib/libdaxctl.so -> libdaxctl.so.1.3.0
lrwxrwxrwx	root/root	usr/lib/libdaxctl.so.1 -> libdaxctl.so.1.3.0
-rwxr-xr-x	root/root	usr/lib/libdaxctl.so.1.3.0
-rwxr-xr-x	root/root	usr/lib/libndctl.la
lrwxrwxrwx	root/root	usr/lib/libndctl.so -> libndctl.so.6.15.0
lrwxrwxrwx	root/root	usr/lib/libndctl.so.6 -> libndctl.so.6.15.0
-rwxr-xr-x	root/root	usr/lib/libndctl.so.6.15.0
drwxr-xr-x	root/root	usr/lib/modprobe.d/
-rw-r--r--	root/root	usr/lib/modprobe.d/nvdimm-security.conf
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libdaxctl.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libndctl.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/ndctl-monitor.service
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/bash-completion/
drwxr-xr-x	root/root	usr/share/bash-completion/completions/
-rw-r--r--	root/root	usr/share/bash-completion/completions/ndctl
drwxr-xr-x	root/root	usr/share/daxctl/
-rw-r--r--	root/root	usr/share/daxctl/daxctl.conf
