# ts/build/packages/netfiles/etc/systemd/system/net.profile-setup.service
[Unit]
Description=ThinStation Profile Setup
Wants=network-online.target
After=network-files.service
Before=shutdown.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/profile-setup
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
