# ts/5.0/packages/base/etc/udev/rules.d/50-modeset.rules
SUBSYSTEM=="module", ACTION=="add", DEVPATH=="/module/i915" RUN+="/etc/udev/scripts/modeset.sh"
SUBSYSTEM=="module", ACTION=="add", DEVPATH=="/module/vmwgfx" RUN+="/etc/udev/scripts/modeset.sh"
SUBSYSTEM=="module", ACTION=="add", DEVPATH=="/module/radeon" RUN+="/etc/udev/scripts/modeset.sh"
SUBSYSTEM=="module", ACTION=="add", DEVPATH=="/module/nouveau" RUN+="/etc/udev/scripts/modeset.sh"

