# ts/5.0/packages/base/etc/udev/rules.d/51-firmware.rules
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/etc/udev/scripts/firmwarelog.sh $env{FIRMWARE}"
