# ts/ports/opt/alsa-utils/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/alsa
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/aconnect
-rwxr-xr-x	root/root	usr/bin/alsabat
-rwxr-xr-x	root/root	usr/bin/alsaloop
-rwxr-xr-x	root/root	usr/bin/alsamixer
-rwxr-xr-x	root/root	usr/bin/alsatplg
-rwxr-xr-x	root/root	usr/bin/alsaucm
-rwxr-xr-x	root/root	usr/bin/amidi
-rwxr-xr-x	root/root	usr/bin/amixer
-rwxr-xr-x	root/root	usr/bin/aplay
-rwxr-xr-x	root/root	usr/bin/aplaymidi
lrwxrwxrwx	root/root	usr/bin/arecord -> aplay
-rwxr-xr-x	root/root	usr/bin/arecordmidi
-rwxr-xr-x	root/root	usr/bin/aseqdump
-rwxr-xr-x	root/root	usr/bin/aseqnet
-rwxr-xr-x	root/root	usr/bin/axfer
-rwxr-xr-x	root/root	usr/bin/iecset
-rwxr-xr-x	root/root	usr/bin/speaker-test
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/alsa-restore.service
-rw-r--r--	root/root	usr/lib/systemd/system/alsa-state.service
drwxr-xr-x	root/root	usr/lib/systemd/system/sound.target.wants/
lrwxrwxrwx	root/root	usr/lib/systemd/system/sound.target.wants/alsa-restore.service -> ../alsa-restore.service
lrwxrwxrwx	root/root	usr/lib/systemd/system/sound.target.wants/alsa-state.service -> ../alsa-state.service
drwxr-xr-x	root/root	usr/lib/udev/
drwxr-xr-x	root/root	usr/lib/udev/rules.d/
-rw-r--r--	root/root	usr/lib/udev/rules.d/89-alsa-ucm.rules
-rw-r--r--	root/root	usr/lib/udev/rules.d/90-alsa-restore.rules
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/fr/
drwxr-xr-x	root/root	usr/man/fr/man8/
-rw-r--r--	root/root	usr/man/fr/man8/alsaconf.8
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/aconnect.1.gz
-rw-r--r--	root/root	usr/man/man1/alsa-info.sh.1.gz
-rw-r--r--	root/root	usr/man/man1/alsabat.1.gz
-rw-r--r--	root/root	usr/man/man1/alsactl.1.gz
-rw-r--r--	root/root	usr/man/man1/alsaloop.1.gz
-rw-r--r--	root/root	usr/man/man1/alsamixer.1.gz
-rw-r--r--	root/root	usr/man/man1/amidi.1.gz
-rw-r--r--	root/root	usr/man/man1/amixer.1.gz
-rw-r--r--	root/root	usr/man/man1/aplay.1.gz
-rw-r--r--	root/root	usr/man/man1/aplaymidi.1.gz
lrwxrwxrwx	root/root	usr/man/man1/arecord.1.gz -> aplay.1.gz
-rw-r--r--	root/root	usr/man/man1/arecordmidi.1.gz
-rw-r--r--	root/root	usr/man/man1/aseqdump.1.gz
-rw-r--r--	root/root	usr/man/man1/aseqnet.1.gz
-rw-r--r--	root/root	usr/man/man1/axfer-list.1.gz
-rw-r--r--	root/root	usr/man/man1/axfer-transfer.1.gz
-rw-r--r--	root/root	usr/man/man1/axfer.1.gz
-rw-r--r--	root/root	usr/man/man1/iecset.1.gz
-rw-r--r--	root/root	usr/man/man1/speaker-test.1.gz
drwxr-xr-x	root/root	usr/man/man7/
-rw-r--r--	root/root	usr/man/man7/alsactl_init.7.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/alsaconf.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/alsa-info.sh
-rwxr-xr-x	root/root	usr/sbin/alsabat-test.sh
-rwxr-xr-x	root/root	usr/sbin/alsaconf
-rwxr-xr-x	root/root	usr/sbin/alsactl
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/alsa/
drwxr-xr-x	root/root	usr/share/alsa/init/
-rw-r--r--	root/root	usr/share/alsa/init/00main
-rw-r--r--	root/root	usr/share/alsa/init/ca0106
-rw-r--r--	root/root	usr/share/alsa/init/default
-rw-r--r--	root/root	usr/share/alsa/init/hda
-rw-r--r--	root/root	usr/share/alsa/init/help
-rw-r--r--	root/root	usr/share/alsa/init/info
-rw-r--r--	root/root	usr/share/alsa/init/test
drwxr-xr-x	root/root	usr/share/alsa/speaker-test/
-rw-r--r--	root/root	usr/share/alsa/speaker-test/sample_map.csv
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/de/
drwxr-xr-x	root/root	usr/share/locale/de/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/de/LC_MESSAGES/alsa-utils.mo
drwxr-xr-x	root/root	usr/share/locale/fr/
drwxr-xr-x	root/root	usr/share/locale/fr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fr/LC_MESSAGES/alsa-utils.mo
drwxr-xr-x	root/root	usr/share/locale/ja/
drwxr-xr-x	root/root	usr/share/locale/ja/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ja/LC_MESSAGES/alsa-utils.mo
-rw-r--r--	root/root	usr/share/locale/ja/LC_MESSAGES/alsaconf.mo
drwxr-xr-x	root/root	usr/share/locale/ru/
drwxr-xr-x	root/root	usr/share/locale/ru/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ru/LC_MESSAGES/alsaconf.mo
drwxr-xr-x	root/root	usr/share/sounds/
drwxr-xr-x	root/root	usr/share/sounds/alsa/
-rw-r--r--	root/root	usr/share/sounds/alsa/Front_Center.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Front_Left.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Front_Right.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Noise.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Rear_Center.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Rear_Left.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Rear_Right.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Side_Left.wav
-rw-r--r--	root/root	usr/share/sounds/alsa/Side_Right.wav
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/alsa/
