# ts/ports/contrib/lm_sensors/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/conf.d/
-rw-r--r--	root/root	etc/conf.d/sensord
-rw-r--r--	root/root	etc/healthd.conf
drwxr-xr-x	root/root	etc/sensors.d/
-rw-r--r--	root/root	etc/sensors3.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/healthd
-rwxr-xr-x	root/root	usr/bin/sensors
-rwxr-xr-x	root/root	usr/bin/sensors-conf-convert
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/sensors/
-rw-r--r--	root/root	usr/include/sensors/error.h
-rw-r--r--	root/root	usr/include/sensors/sensors.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libsensors.a
lrwxrwxrwx	root/root	usr/lib/libsensors.so -> libsensors.so.5
lrwxrwxrwx	root/root	usr/lib/libsensors.so.5 -> libsensors.so.5.0.0
-rwxr-xr-x	root/root	usr/lib/libsensors.so.5.0.0
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/fancontrol.service
-rw-r--r--	root/root	usr/lib/systemd/system/healthd.service
-rw-r--r--	root/root	usr/lib/systemd/system/lm_sensors.service
-rw-r--r--	root/root	usr/lib/systemd/system/sensord.service
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/fancontrol
-rwxr-xr-x	root/root	usr/sbin/isadump
-rwxr-xr-x	root/root	usr/sbin/isaset
-rwxr-xr-x	root/root	usr/sbin/pwmconfig
-rwxr-xr-x	root/root	usr/sbin/sensors-detect
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/sensors.1.gz
drwxr-xr-x	root/root	usr/share/man/man3/
-rw-r--r--	root/root	usr/share/man/man3/libsensors.3.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/sensors.conf.5.gz
lrwxrwxrwx	root/root	usr/share/man/man5/sensors3.conf.5.gz -> sensors.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/fancontrol.8.gz
-rw-r--r--	root/root	usr/share/man/man8/isadump.8.gz
-rw-r--r--	root/root	usr/share/man/man8/isaset.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pwmconfig.8.gz
-rw-r--r--	root/root	usr/share/man/man8/sensors-conf-convert.8.gz
-rw-r--r--	root/root	usr/share/man/man8/sensors-detect.8.gz
