# ts/ports/components/opencryptoki/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/ld.so.conf.d/
-rw-r--r--	root/root	etc/ld.so.conf.d/opencryptoki-i686.conf
drwxr-xr-x	root/root	etc/rc.d/
drwxr-xr-x	root/root	etc/rc.d/init.d/
-rwxr-xr-x	root/root	etc/rc.d/init.d/pkcsslotd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/opencryptoki/
-rw-r--r--	root/root	usr/include/opencryptoki/apiclient.h
-rw-r--r--	root/root	usr/include/opencryptoki/pkcs11.h
-rw-r--r--	root/root	usr/include/opencryptoki/pkcs11types.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/opencryptoki/
lrwxrwxrwx	root/root	usr/lib/opencryptoki/PKCS11_API.so -> libopencryptoki.so
-rwxr-xr-x	root/root	usr/lib/opencryptoki/libopencryptoki.la
lrwxrwxrwx	root/root	usr/lib/opencryptoki/libopencryptoki.so -> libopencryptoki.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/opencryptoki/libopencryptoki.so.0 -> libopencryptoki.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/opencryptoki/libopencryptoki.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/opencryptoki/methods -> /usr/sbin
drwxr-xr-x	root/root	usr/lib/opencryptoki/stdll/
lrwxrwxrwx	root/root	usr/lib/opencryptoki/stdll/PKCS11_CCA.so -> libpkcs11_cca.so
lrwxrwxrwx	root/root	usr/lib/opencryptoki/stdll/PKCS11_SW.so -> libpkcs11_sw.so
-rwxr-xr-x	root/root	usr/lib/opencryptoki/stdll/libpkcs11_cca.la
lrwxrwxrwx	root/root	usr/lib/opencryptoki/stdll/libpkcs11_cca.so -> libpkcs11_cca.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/opencryptoki/stdll/libpkcs11_cca.so.0 -> libpkcs11_cca.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/opencryptoki/stdll/libpkcs11_cca.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/opencryptoki/stdll/libpkcs11_sw.la
lrwxrwxrwx	root/root	usr/lib/opencryptoki/stdll/libpkcs11_sw.so -> libpkcs11_sw.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/opencryptoki/stdll/libpkcs11_sw.so.0 -> libpkcs11_sw.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/opencryptoki/stdll/libpkcs11_sw.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkcs11/
lrwxrwxrwx	root/root	usr/lib/pkcs11/PKCS11_API.so -> ../opencryptoki/libopencryptoki.so
lrwxrwxrwx	root/root	usr/lib/pkcs11/libopencryptoki.so -> ../opencryptoki/libopencryptoki.so
lrwxrwxrwx	root/root	usr/lib/pkcs11/methods -> /usr/sbin
lrwxrwxrwx	root/root	usr/lib/pkcs11/stdll -> ../opencryptoki/stdll/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/pkcs11_startup.1.gz
-rw-r--r--	root/root	usr/man/man1/pkcs_slot.1.gz
-rw-r--r--	root/root	usr/man/man1/pkcsconf.1.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/pk_config_data.5.gz
drwxr-xr-x	root/root	usr/man/man7/
-rw-r--r--	root/root	usr/man/man7/opencryptoki.7.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/pkcsslotd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/pkcs11_startup
-rwxr-xr-x	root/root	usr/sbin/pkcs_slot
-rwxr-xr-x	root/root	usr/sbin/pkcsconf
-rwxr-xr-x	root/root	usr/sbin/pkcsslotd
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/opencryptoki/
drwxr-xr-x	root/root	var/lib/opencryptoki/ccatok/
drwxr-xr-x	root/root	var/lib/opencryptoki/swtok/
drwxr-xr-x	root/root	var/lock/
drwxrwx---	root/pkcs11	var/lock/opencryptoki/
