# ts/build/packages/base/etc/systemd/user/homedir.service
[Unit]
Description=Create User Home Directory

[Service]
ExecStart=/etc/systemd/user/session open

[Install]
WantedBy=default.target
