# ts/build/packages/alsa/build/extra/lib/systemd/system/audio.init.service
[Unit]
Description=Configure Audio Settings
DefaultDependencies=no
Conflicts=shutdown.target
After=profile_setup.service net.profle_setup.service
After=alsa-restore.service
Before=session-setup.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/lib/thinstation/system/audio
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
