# ts/5.1/packages/pulseaudio/build/extra/etc/xdg/autostart/pulseaudio.desktop
[Desktop Entry]
Type=Application
Name=Pulse Audio Setup
Exec=/bin/start-pulse-x11
OnlyShowIn=GNOME;
NoDisplay=true
X-GNOME-Autostart-Phase=Initialization
X-GNOME-Autostart-Notify=true
X-GNOME-AutoRestart=true
