# ts/build/packages/xorg7/lib/systemd/system/xconf.service
[Unit]
Description=Configure Thinstation Xorg Settings
DefaultDependencies=no
Conflicts=shutdown.target
After=session.service
Before=session-setup.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/lib/thinstation/system/xconf
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
