# ts/ports/components/geoclue/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.GeoClue2.Agent.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.GeoClue2.conf
drwxr-xr-x	root/root	etc/geoclue/
-rw-r--r--	root/root	etc/geoclue/geoclue.conf
drwxr-xr-x	root/root	etc/xdg/
drwxr-xr-x	root/root	etc/xdg/autostart/
-rw-r--r--	root/root	etc/xdg/autostart/geoclue-demo-agent.desktop
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/libgeoclue-2.0/
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-client.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-enum-types.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-enums.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-helpers.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-location.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-manager.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-simple.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/geoclue.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/girepository-1.0/
-rw-r--r--	root/root	usr/lib/girepository-1.0/Geoclue-2.0.typelib
lrwxrwxrwx	root/root	usr/lib/libgeoclue-2.so -> libgeoclue-2.so.0
lrwxrwxrwx	root/root	usr/lib/libgeoclue-2.so.0 -> libgeoclue-2.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libgeoclue-2.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/geoclue-2.0.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libgeoclue-2.0.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/geoclue.service
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/geoclue
drwxr-xr-x	root/root	usr/libexec/geoclue-2.0/
drwxr-xr-x	root/root	usr/libexec/geoclue-2.0/demos/
-rwxr-xr-x	root/root	usr/libexec/geoclue-2.0/demos/agent
-rwxr-xr-x	root/root	usr/libexec/geoclue-2.0/demos/where-am-i
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/applications/
-rw-r--r--	root/root	usr/share/applications/geoclue-demo-agent.desktop
-rw-r--r--	root/root	usr/share/applications/geoclue-where-am-i.desktop
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/interfaces/
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Agent.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Client.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Location.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Manager.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.xml
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.GeoClue2.service
drwxr-xr-x	root/root	usr/share/gir-1.0/
-rw-r--r--	root/root	usr/share/gir-1.0/Geoclue-2.0.gir
drwxr-xr-x	root/root	usr/share/gtk-doc/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/geoclue/
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/gdbus-org.freedesktop.GeoClue2.Agent.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/gdbus-org.freedesktop.GeoClue2.Client.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/gdbus-org.freedesktop.GeoClue2.Location.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/gdbus-org.freedesktop.GeoClue2.Manager.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/geoclue-gclue-enums.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/geoclue.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/ix01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/left-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/license.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/ref-agent-dbus.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/ref-dbus.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/right-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/up-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/geoclue/up.png
drwxr-xr-x	root/root	usr/share/gtk-doc/html/libgeoclue/
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/GClueClient.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/GClueClientProxy.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/GClueLocation.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/GClueLocationProxy.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/GClueManager.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/GClueManagerProxy.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/GClueSimple.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/annotation-glossary.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/ch01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/ch02.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/ch03.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/ch04.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/ch05.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/ix01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/left-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/libgeoclue-gclue-enums.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/libgeoclue.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/license.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/right-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/up-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libgeoclue/up.png
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/geoclue.5.gz
drwxr-xr-x	root/root	usr/share/vala/
drwxr-xr-x	root/root	usr/share/vala/vapi/
-rw-r--r--	root/root	usr/share/vala/vapi/libgeoclue-2.0.deps
-rw-r--r--	root/root	usr/share/vala/vapi/libgeoclue-2.0.vapi
