# ts/ports/components/geoclue/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.GeoClue2.Agent.conf
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.GeoClue2.conf
drwxr-xr-x	root/root	etc/geoclue/
-rw-r--r--	root/root	etc/geoclue/geoclue.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/libgeoclue-2.0/
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-client.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-enum-types.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-enums.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-helpers.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-location.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-manager.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/gclue-simple.h
-rw-r--r--	root/root	usr/include/libgeoclue-2.0/geoclue.h
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/geoclue
drwxr-xr-x	root/root	usr/lib/geoclue-2.0/
drwxr-xr-x	root/root	usr/lib/geoclue-2.0/demos/
-rwxr-xr-x	root/root	usr/lib/geoclue-2.0/demos/where-am-i
drwxr-xr-x	root/root	usr/lib/girepository-1.0/
-rw-r--r--	root/root	usr/lib/girepository-1.0/Geoclue-2.0.typelib
-rwxr-xr-x	root/root	usr/lib/libgeoclue-2.la
lrwxrwxrwx	root/root	usr/lib/libgeoclue-2.so -> libgeoclue-2.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libgeoclue-2.so.0 -> libgeoclue-2.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libgeoclue-2.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/geoclue-2.0.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libgeoclue-2.0.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/geoclue.service
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/applications/
-rw-r--r--	root/root	usr/share/applications/geoclue-demo-agent.desktop
-rw-r--r--	root/root	usr/share/applications/geoclue-where-am-i.desktop
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/interfaces/
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Agent.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Client.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Location.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Manager.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.xml
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.GeoClue2.service
drwxr-xr-x	root/root	usr/share/gir-1.0/
-rw-r--r--	root/root	usr/share/gir-1.0/Geoclue-2.0.gir
