# ts/build/packages/vmview/build/install
#!/bin/bash

gzname=`basename $1`
INSTALLDIR=./packages/vmview

components="\
--set-setting vmware-horizon-smartcard smartcardEnable yes \
--set-setting vmware-horizon-rtav rtavEnable yes \
--set-setting vmware-horizon-usb usbEnable yes \
--set-setting vmware-horizon-virtual-printing tpEnable yes \
--set-setting vmware-horizon-tsdr tsdrEnable yes \
--set-setting vmware-horizon-mmr mmrEnable yes \
--set-setting vmware-horizon-scannerclient scannerEnable yes \
--set-setting vmware-horizon-media-provider mediaproviderEnable no"

new_file()
{
        if [ -n "`stat -c '%Y' $1 | grep -v \"$bday\"`" ]; then
                return 0
        else
                return 1
        fi
}

etcs="passwd group nsswitch.conf issue ld.so.conf hosts"
bins="stat attr awk basename busybox bash cat chgrp chmod chown cp cut date dd df \
	dirname echo ed egrep expr false find gawk grep head id init od \
	install ip ln locale ls mkdir mktemp more mount mv pidof ping \
	ps red rm rmdir sed sh sleep sort ss tar tc telinit touch tr \
	true uname vi wc whoami ldd file ldconfig strace getconf"
dirs="	proc \
	sys \
	etc/udev/rules.d \
	etc/init.d \
	bin \
	lib \
	tmp \
	dev/bus/usb \
	usr/local \
	sbin \
        lib/applications \
	etc/init.d \
	etc/rc0.d \
	etc/rc1.d \
	etc/rc2.d \
	etc/rc3.d \
	etc/rc4.d \
	etc/rc5.d \
	etc/rc6.d \
	build "
clean_dirs="	dev \
		proc \
		var \
		tmp \
		sys \
		lib/python2.7 \
		lib/vmware-installer-horizon \
		etc/vmware-installer-horizon \
		etc/ld.so.cache \
		lib/vmware/gcc \
		lib/doc/vmware-horizon-client/debug \
		lib/doc/vmware-horizon-client/patches \
		etc/rc0.d \
		etc/rc1.d \
		etc/rc2.d \
		etc/rc3.d \
		etc/rc4.d \
		etc/rc5.d \
		etc/rc6.d "
libs="	/usr/lib/libffi.so.6 \
	/usr/lib/libstdc++.so.6 \
	/usr/lib/libFLAC.so.8 \
	/usr/lib/libXau.so.6 \
	/usr/lib/libXdmcp.so.6 \
	/usr/lib/libXft.so.2 \
	/usr/lib/libogg.so.0 \
	/usr/lib/libvorbis.so.0 \
	/usr/lib/libvorbisenc.so.2 \
	/usr/lib/libatspi.so.0 \
	/usr/lib/libdrm.so.2 \
	/usr/lib/libfribidi.so.0 \
	/usr/lib/libharfbuzz.so.0 \
	/lib/liblzo2.so.2 \
	/usr/lib/liborc-0.4.so.0 \
	/usr/lib/libpng16.so.16 \
	/usr/lib/libsndfile.so.1 \
	/usr/lib/libwayland-client.so.0 \
	/lib/libnss_files.so.2 \
	/lib/libutil.so.1 \
	/usr/lib/libsqlite3.so.0 \
	/usr/lib/libgssapi_krb5.so.2 \
	/lib/libkrb5.so.3 \
	/lib/libk5crypto.so.3 \
	/lib/libcom_err.so.2 \
	/usr/lib/libgnutls.so.30 \
	/usr/lib/libavahi-common.so.3 \
	/usr/lib/libavahi-client.so.3 \
	/lib/libkrb5support.so.0 \
	/lib/libkeyutils.so.1 \
	/lib/libresolv.so.2 \
	/usr/lib/libnettle.so.8 \
	/usr/lib/libhogweed.so.6 \
	/usr/lib/libdbus-1.so.3 \
	/usr/lib/libsystemd.so.0 \
	/usr/lib/libgcc_s.so.1"

# Create all the directories that we are going to need
for repack in python cups expat sane libgphoto2 v4l-utils libtiff libjpeg-turbo libxkbcommon libepoxy \
	libusb libwebp zstd libtool libexif systemd libgcrypt libgpg-error gtk gtk3 wayland \
	openssl gst-plugins-base gstreamer pygobject-3 \
	dbus-python gdb meson pip2 scons setuptools2 netaddr pyxdg sqlalchemy urllib3 \
	libpng12 pcsc-lite atkmm cairomm ffmpeg libsigc++ pangomm util-linux at-spi2-atk \
	atk gdk-pixbuf glib glibmm gtk gtkmm libxml2 pango alsa-lib cairo fontconfig freetype libva \
	nspr nss pulseaudio libvdpau xorg-libpixman xorg-libx11 xorg-libxcb \
	xorg-libxcomposite xorg-libxcursor xorg-libxdamage xorg-libxext xorg-libxfixes xorg-libxi \
	xorg-libxinerama xorg-libxkbfile xorg-libxrandr xorg-libxrender xorg-libxscrnsaver xorg-libxtst ; do
	repackage $repack -t $INSTALLDIR
	cp -a $INSTALLDIR/$repack/* $INSTALLDIR/.
	rm -rf $INSTALLDIR/$repack
done

for mdir in $dirs; do
        mkdir -p $INSTALLDIR/$mdir
done

# Fill in the bin folder
for cbin in $bins ; do
        cp `which $cbin` $INSTALLDIR/bin/.
	if [ -n "`file $INSTALLDIR/bin/$cbin |grep -oe ELF`" ]; then
		libs="$libs \n`ldd $INSTALLDIR/bin/$cbin |grep -v 'linux-vdso' |grep -v 'not a' |cut -d '>' -f2 |cut -d '(' -f1 |sed -E s/[[:space:]]//g`"
	fi
done
cp /usr/bin/which $INSTALLDIR/bin/.

ln -sf python2 $INSTALLDIR/bin/python

if [ -e /bin/uname-auto ] ; then
        cp /bin/uname-auto $INSTALLDIR/bin/.
fi
echo -e "#!/bin/sh\nexit 0" > $INSTALLDIR/bin/install_initd
echo -e "#!/bin/sh\nexit 0" > $INSTALLDIR/bin/update-rc.d
echo -e "#!/bin/sh\nexit 0" > $INSTALLDIR/bin/chkconfig
chmod +x $INSTALLDIR/bin/install_initd
chmod +x $INSTALLDIR/bin/update-rc.d
chmod +x $INSTALLDIR/bin/chkconfig
touch $INSTALLDIR/bin/firefox
touch $INSTALLDIR/bin/chrome

# Fill in the lib folder
for clib in `echo -e $libs |sort -u`; do
        cp $clib $INSTALLDIR/lib/.
done

# Fill in the etc folder
for cetc in $etcs ; do
        cp /etc/$cetc $INSTALLDIR/etc/.
done
cat /proc/mounts |grep -Ev 'gvfs|mapper|media' > $INSTALLDIR/etc/mtab

# Add some links to fool any tests
ln -sf /sbin $INSTALLDIR/usr/local/sbin
ln -sf /bin $INSTALLDIR/usr/local/bin
ln -sf /lib $INSTALLDIR/usr/local/lib
ln -sf /sbin $INSTALLDIR/usr/sbin
ln -sf /bin $INSTALLDIR/usr/bin
ln -sf /lib $INSTALLDIR/usr/lib
ln -sf /lib $INSTALLDIR/usr/share
ln -sf /lib $INSTALLDIR/lib64
ln -sf /lib $INSTALLDIR/usr/lib64

# Pop busybox and make installer executable
chroot $INSTALLDIR /bin/busybox --install -s
chroot $INSTALLDIR /bin/ldconfig
cp /downloads/$gzname $INSTALLDIR/.
chmod +x $INSTALLDIR/$gzname

# This will mark all the files before we run install, so that we know what to remove afterwards.
timestamp="`date +%Y%m%d`0000.00"
for file in `find $INSTALLDIR`; do
	touch -c -h -t $timestamp $file
done
bday=`stat -c %Y $INSTALLDIR/$gzname`

# Start the installer in a chroot session.
mount -t proc proc $INSTALLDIR/proc
mount -t sysfs sys $INSTALLDIR/sys
mount -t devtmpfs dev $INSTALLDIR/dev

# Run the installer
chroot $INSTALLDIR strace -f -o /debug ./$gzname -I --stop-services --eulas-agreed --console --required $components
chroot $INSTALLDIR ./$gzname -t -I

for mount in dev proc sys; do
	while mountpoint $INSTALLDIR/$mount >/dev/null 2>&1 ; do
		umount $INSTALLDIR/$mount
		sleep 1
	done
done

cp $INSTALLDIR/tmp/vmware-root $INSTALLDIR/build/vmware-log -a
cp $INSTALLDIR/var/log/vmware/* $INSTALLDIR/build/vmware-log/.
cp $INSTALLDIR/debug $INSTALLDIR/build/vmware-log/install.strace

# Remove folders and files that we know must go.
for dir in $clean_dirs ; do
        rm -rf $INSTALLDIR/$dir
done
rm $INSTALLDIR/$gzname
rm $INSTALLDIR/lib64
rm $INSTALLDIR/debug

# Pin open empty folders, they may get used on a running system
for dir in `find $INSTALLDIR -type d -empty`; do
        touch $dir/.gitignore
done

# Remove any files not placed by the installer according to the bday set earlier
for dir in `find $INSTALLDIR -mindepth 1 -maxdepth 1 -type d ! -name build`; do
	for file in `find $dir -type f -o -type l`; do
		if ! new_file $file; then
			rm $file
		fi
	done
done

# Merge static files and folders
cp -Prf $INSTALLDIR/build/extra/* $INSTALLDIR/.

# Fix links, relative linking is bad.
for link in `find $INSTALLDIR/etc/init.d -type l`; do
	dest=`readlink $link`
	newdest=`echo $dest |sed -e "s|[.]\+/[.]\+/usr/lib64|/lib|g"`
	rm $link
	ln -sf $newdest $link
done
for link in `find $INSTALLDIR/bin -type l`; do
        dest=`readlink $link`
        newdest=`echo $dest |sed -e "s|[.]\+/lib64|/lib|g"`
        rm $link
        ln -sf $newdest $link
done
for link in `find $INSTALLDIR/sbin -type l`; do
        dest=`readlink $link`
        newdest=`echo $dest |sed -e "s|[.]\+/usr/lib64|/lib|g"`
        rm $link
        ln -sf $newdest $link
done

# Delete libraries we already have in our build system
#for lib in `find $INSTALLDIR ! -type d -name lib\*|grep -v glibmm`; do
#	if readlink $lib >/dev/null; then
#		if [ -z $(find-dep-files $(basename $lib) $INSTALLDIR 2>/dev/null) ]; then
#			rm $lib
#		fi
#	elif prt-get fsearch `basename $lib` > /dev/null; then
#		rm $lib
#	fi
#done

# Fix typos and work around missing features in busybox
if [ -e $INSTALLDIR/lib/vmware/view/env/env_utils.sh ]; then
  sed -i -e 's/glibcxx=.*/glibcxx=1/' $INSTALLDIR/lib/vmware/view/env/env_utils.sh
else
  sed -i -e 's/$ENABLE_SEAMLESS_WINDOW = /"$ENABLE_SEAMLESS_WINDOW" = /' \
	-e 's/glibcxx=.*/glibcxx=1/' $INSTALLDIR/bin/vmware-view
fi

# Remove broken link
rm $INSTALLDIR/lib/vmware/libudev.so.0

# Remove any now empty folders
find $INSTALLDIR -type d -empty -delete

# Remove the debugging symbols ## VMWare ??
for elf in `find $INSTALLDIR -type f |xargs file |grep -e ELF |cut -d ":" -f1` ; do
	strip --strip-all $elf
done

ln -s ../lib/doc/vmware-horizon-client/VMware-Horizon-Client-EULA-en.txt $INSTALLDIR/build/license
