# ts/build/packages/ntpdate/etc/systemd/system/ntpdate.service
[Unit]
Description=Thistation ntpd
After=fastboot.service NetworkManager-wait-online.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/ntpdate start
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
