# ts/build/packages/vmviewpcoip/build/extra/lib/systemd/system/vmware.view.usbd.service
[Unit]
Description=VMware Horizon Client USB Sharing Service
Requires=vmware.usbarbitrator.service
After=vmware.usbarbitrator.service

[Service]
ExecStart=/usr/lib/vmware/view/usb/vmware-view-usbd -f

[Install]
WantedBy=multi-user.target
