# ts/build/packages/vmviewpcoip/build/extra/lib/systemd/system/vmware.usbarbitrator.service
[Unit]
Description=VMware Horizon Client USB Arbitrator

[Service]
ExecStart=/usr/lib/vmware/view/usb/vmware-usbarbitrator -f
ExecKill=/usr/lib/vmware/view/usb/vmware-usbarbitrator --kill

[Install]
WantedBy=multi-user.target
