# ts/ports/components/telnet-bsd/telnet.socket
[Unit]
Description=Telnet per-connection socket

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
