# ts/ports/components/telnet-bsd/telnet@.service
[Unit]
Description=Telnet Per-Connection Server

[Service]
ExecStart=-/usr/bin/in.telnetd
User=root
StandardInput=socket
