# ts/build/packages/netfiles/lib/systemd/system/network_files.service
[Unit]
Description=Thistation Network Files
Wants=network-online.target
After=profile_setup.service network-online.target
Before=net.profile-setup.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/lib/thinstation/system/network_files
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
