# ts/ports/components/xinetd-TS/.footprint
drwxr-xr-x	root/root	ts/
drwxr-xr-x	root/root	ts/components/
drwxr-xr-x	root/root	ts/components/xinetd/
drwxr-xr-x	root/root	ts/components/xinetd/etc/
drwxr-xr-x	root/root	ts/components/xinetd/etc/rc.d/
-rw-rw-r--	root/root	ts/components/xinetd/etc/xinetd.conf.sample
drwxr-xr-x	root/root	ts/components/xinetd/sbin/
-rwxr-xr-x	root/root	ts/components/xinetd/sbin/itox
-rwxr-xr-x	root/root	ts/components/xinetd/sbin/xconv.pl
-rwxr-xr-x	root/root	ts/components/xinetd/sbin/xinetd
drwxr-xr-x	root/root	ts/components/xinetd/var/
drwxr-xr-x	root/root	ts/components/xinetd/var/log/
-rw-------	root/root	ts/components/xinetd/var/log/xinetd.log (EMPTY)
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/xinetd.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/itox.8.gz
-rw-r--r--	root/root	usr/man/man8/xconv.pl.8.gz
-rw-r--r--	root/root	usr/man/man8/xinetd.8.gz
-rw-r--r--	root/root	usr/man/man8/xinetd.log.8.gz
