# ts/build/packages/lprng/etc/systemd/system/lprng.service
[Unit]
Description=ThinStation lprng
After=profile-setup.service pkg.service
Before=display-manager.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/lprng init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
