# ts/ports/components/acpid-TS/.footprint
drwxr-sr-x	root/4135	ts/
drwxr-sr-x	root/4135	ts/components/
drwxr-sr-x	root/4135	ts/components/acpid/
drwxr-sr-x	root/4135	ts/components/acpid/bin/
-rwxr-xr-x	root/4135	ts/components/acpid/bin/acpi_listen
-rwxr-xr-x	root/4135	ts/components/acpid/bin/acpid
-rwxr-xr-x	root/4135	ts/components/acpid/bin/kacpimon
drwxr-sr-x	root/4135	ts/components/acpid/etc/
drwxr-sr-x	root/4135	ts/components/acpid/etc/acpi/
-rwxr-xr-x	root/4135	ts/components/acpid/etc/acpi/default.sh
drwxr-sr-x	root/4135	ts/components/acpid/etc/acpi/events/
-rw-r--r--	root/4135	ts/components/acpid/etc/acpi/events/default
drwxr-sr-x	root/4135	ts/components/acpid/etc/rc.d/
-rwxr-xr-x	root/4135	ts/components/acpid/etc/rc.d/acpid
