# ts/ports/components/v4l-utils/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc_keymaps/
drwxr-xr-x	root/root	etc/rc_keymaps/protocols/
-rw-r--r--	root/root	etc/rc_maps.cfg
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/cec-compliance
-rwxr-xr-x	root/root	usr/bin/cec-ctl
-rwxr-xr-x	root/root	usr/bin/cec-follower
-rwxr-xr-x	root/root	usr/bin/cx18-ctl
-rwxr-xr-x	root/root	usr/bin/decode_tm6000
-rwxr-xr-x	root/root	usr/bin/dvb-fe-tool
-rwxr-xr-x	root/root	usr/bin/dvb-format-convert
-rwxr-xr-x	root/root	usr/bin/dvbv5-daemon
-rwxr-xr-x	root/root	usr/bin/dvbv5-scan
-rwxr-xr-x	root/root	usr/bin/dvbv5-zap
-rwxr-xr-x	root/root	usr/bin/ir-ctl
-rwxr-xr-x	root/root	usr/bin/ir-keytable
-rwxr-xr-x	root/root	usr/bin/ivtv-ctl
-rwxr-xr-x	root/root	usr/bin/media-ctl
-rwxr-xr-x	root/root	usr/bin/qvidcap
-rwxr-xr-x	root/root	usr/bin/rds-ctl
-rwxr-xr-x	root/root	usr/bin/v4l2-compliance
-rwxr-xr-x	root/root	usr/bin/v4l2-ctl
-rwxr-xr-x	root/root	usr/bin/v4l2-sysfs-path
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/libdvbv5/
-rw-r--r--	root/root	usr/include/libdvbv5/atsc_eit.h
-rw-r--r--	root/root	usr/include/libdvbv5/atsc_header.h
-rw-r--r--	root/root	usr/include/libdvbv5/cat.h
-rw-r--r--	root/root	usr/include/libdvbv5/countries.h
-rw-r--r--	root/root	usr/include/libdvbv5/crc32.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_atsc_service_location.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_ca.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_ca_identifier.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_cable_delivery.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_event_extended.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_event_short.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_extension.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_frequency_list.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_hierarchy.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_isdbt_delivery.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_language.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_logical_channel.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_network_name.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_partial_reception.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_sat.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_service.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_t2_delivery.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_terrestrial_delivery.h
-rw-r--r--	root/root	usr/include/libdvbv5/desc_ts_info.h
-rw-r--r--	root/root	usr/include/libdvbv5/descriptors.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-demux.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-dev.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-fe.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-file.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-frontend.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-log.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-sat.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-scan.h
-rw-r--r--	root/root	usr/include/libdvbv5/dvb-v5-std.h
-rw-r--r--	root/root	usr/include/libdvbv5/eit.h
-rw-r--r--	root/root	usr/include/libdvbv5/header.h
-rw-r--r--	root/root	usr/include/libdvbv5/libdvb-version.h
-rw-r--r--	root/root	usr/include/libdvbv5/mgt.h
-rw-r--r--	root/root	usr/include/libdvbv5/mpeg_es.h
-rw-r--r--	root/root	usr/include/libdvbv5/mpeg_pes.h
-rw-r--r--	root/root	usr/include/libdvbv5/mpeg_ts.h
-rw-r--r--	root/root	usr/include/libdvbv5/nit.h
-rw-r--r--	root/root	usr/include/libdvbv5/pat.h
-rw-r--r--	root/root	usr/include/libdvbv5/pmt.h
-rw-r--r--	root/root	usr/include/libdvbv5/sdt.h
-rw-r--r--	root/root	usr/include/libdvbv5/vct.h
-rw-r--r--	root/root	usr/include/libv4l-plugin.h
-rw-r--r--	root/root	usr/include/libv4l1-videodev.h
-rw-r--r--	root/root	usr/include/libv4l1.h
-rw-r--r--	root/root	usr/include/libv4l2.h
-rw-r--r--	root/root	usr/include/libv4l2rds.h
-rw-r--r--	root/root	usr/include/libv4lconvert.h
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/libdvbv5.la
lrwxrwxrwx	root/root	usr/lib/libdvbv5.so -> libdvbv5.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libdvbv5.so.0 -> libdvbv5.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libdvbv5.so.0.0.0
drwxr-xr-x	root/root	usr/lib/libv4l/
-rwxr-xr-x	root/root	usr/lib/libv4l/ov511-decomp
-rwxr-xr-x	root/root	usr/lib/libv4l/ov518-decomp
drwxr-xr-x	root/root	usr/lib/libv4l/plugins/
-rwxr-xr-x	root/root	usr/lib/libv4l/plugins/libv4l-mplane.la
-rwxr-xr-x	root/root	usr/lib/libv4l/plugins/libv4l-mplane.so
-rwxr-xr-x	root/root	usr/lib/libv4l/v4l1compat.la
-rwxr-xr-x	root/root	usr/lib/libv4l/v4l1compat.so
-rwxr-xr-x	root/root	usr/lib/libv4l/v4l2convert.la
-rwxr-xr-x	root/root	usr/lib/libv4l/v4l2convert.so
-rwxr-xr-x	root/root	usr/lib/libv4l1.la
lrwxrwxrwx	root/root	usr/lib/libv4l1.so -> libv4l1.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libv4l1.so.0 -> libv4l1.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libv4l1.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libv4l2.la
lrwxrwxrwx	root/root	usr/lib/libv4l2.so -> libv4l2.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libv4l2.so.0 -> libv4l2.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libv4l2.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libv4l2rds.la
lrwxrwxrwx	root/root	usr/lib/libv4l2rds.so -> libv4l2rds.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libv4l2rds.so.0 -> libv4l2rds.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libv4l2rds.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libv4lconvert.la
lrwxrwxrwx	root/root	usr/lib/libv4lconvert.so -> libv4lconvert.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libv4lconvert.so.0 -> libv4lconvert.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libv4lconvert.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libdvbv5.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libv4l1.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libv4l2.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libv4l2rds.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libv4lconvert.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
drwxr-xr-x	root/root	usr/lib/systemd/system/systemd-udevd.service.d/
-rw-r--r--	root/root	usr/lib/systemd/system/systemd-udevd.service.d/50-rc_keymap.conf
drwxr-xr-x	root/root	usr/lib/udev/
drwxr-xr-x	root/root	usr/lib/udev/rc_keymaps/
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/adstech_dvb_t_pci.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/af9005.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/alink_dtu_m.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/allwinner_ba10_tv_box.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/allwinner_i12_a20_tv_box.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/anysee.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/apac_viewcomp.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/astrometa_t2hybrid.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/asus_pc39.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/asus_ps3_100.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/ati_tv_wonder_hd_600.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/ati_x10.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avermedia.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avermedia_a16d.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avermedia_cardbus.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avermedia_dvbt.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avermedia_m135a.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avermedia_m733a_rm_k6.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avermedia_rm_ks.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/avertv_303.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/az6027.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/azurewave_ad_tu700.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/beelink_gs1.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/behold.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/behold_columbus.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/budget_ci_old.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/cec.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/cinergy.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/cinergy_1400.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/cinergyt2.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/d680_dmb.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/delock_61959.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dib0700_nec.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dib0700_rc5.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dibusb.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/digitalnow_tinytwin.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/digittrade.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/digitv.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dish_network.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dm1105_nec.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dntv_live_dvb_t.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dntv_live_dvbt_pro.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dtt200u.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dvbsky.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dvico_mce.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/dvico_portable.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/em_terratec.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/encore_enltv.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/encore_enltv2.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/encore_enltv_fm53.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/evga_indtube.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/eztv.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/flydvb.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/flyvideo.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/fusionhdtv_mce.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/gadmei_rm008z.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/geekbox.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/genius_tvgo_a11mce.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/gotview7135.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/haupp.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/hauppauge.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/hisi_poplar.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/hisi_tv_demo.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/imon_mce.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/imon_pad.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/imon_rsc.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/iodata_bctv7e.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/it913x_v1.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/it913x_v2.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/kaiomy.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/khadas.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/kii_pro.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/kworld_315u.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/kworld_pc150u.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/kworld_plus_tv_analog.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/leadtek_y04g0051.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/lme2510.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/manli.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/mce_keyboard.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/medion_x10.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/medion_x10_digitainer.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/medion_x10_or2x.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/megasky.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/msi_digivox_ii.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/msi_digivox_iii.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/msi_tvanywhere.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/msi_tvanywhere_plus.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/nebula.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/nec_terratec_cinergy_xs.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/norwood.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/npgtech.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/odroid.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/opera1.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pctv_sedna.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pinnacle310e.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pinnacle_color.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pinnacle_grey.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pinnacle_pctv_hd.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pixelview.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pixelview_002t.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pixelview_mk12.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pixelview_new.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/powercolor_real_angel.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/proteus_2309.toml
drwxr-xr-x	root/root	usr/lib/udev/rc_keymaps/protocols/
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/grundig.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/imon_rsc.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/manchester.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/pulse_distance.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/pulse_length.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/raw.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/rc_mm.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/samsung36.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/protocols/xbox-dvd.o
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/purpletv.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/pv951.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/rc6_mce.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/real_audio_220_32_keys.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/reddo.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/snapstream_firefly.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/streamzap.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/su3000.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tango.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tanix_tx3mini.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tanix_tx5max.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tbs_nec.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/technisat_ts35.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/technisat_usb2.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/terratec_cinergy_c_pci.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/terratec_cinergy_s2_hd.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/terratec_cinergy_xs.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/terratec_slim.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/terratec_slim_2.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tevii_nec.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tivo.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/total_media_in_hand.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/total_media_in_hand_02.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/trekstor.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tt_1500.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/tvwalkertwin.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/twinhan_dtv_cab_ci.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/twinhan_vp1027_dvbs.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/vega_s9x.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/videomate_k100.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/videomate_s350.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/videomate_tv_pvr.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/vp702x.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/wetek_hub.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/wetek_play2.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/winfast.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/winfast_usbii_deluxe.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/wobo_i5.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/x96max.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/xbox_dvd.toml
-rw-r--r--	root/root	usr/lib/udev/rc_keymaps/zx_irdec.toml
drwxr-xr-x	root/root	usr/lib/udev/rules.d/
-rw-r--r--	root/root	usr/lib/udev/rules.d/70-infrared.rules
lrwxrwxrwx	root/root	usr/lib/v4l1compat.so -> libv4l/v4l1compat.so
lrwxrwxrwx	root/root	usr/lib/v4l2convert.so -> libv4l/v4l2convert.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/v4l2-dbg
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/applications/
-rw-r--r--	root/root	usr/share/applications/qvidcap.desktop
drwxr-xr-x	root/root	usr/share/icons/
drwxr-xr-x	root/root	usr/share/icons/hicolor/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/16x16/apps/qvidcap.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/24x24/apps/qvidcap.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/32x32/apps/qvidcap.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/64x64/
drwxr-xr-x	root/root	usr/share/icons/hicolor/64x64/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/64x64/apps/qvidcap.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/scalable/apps/qvidcap.svg
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/ca/
drwxr-xr-x	root/root	usr/share/locale/ca/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ca/LC_MESSAGES/libdvbv5.mo
-rw-r--r--	root/root	usr/share/locale/ca/LC_MESSAGES/v4l-utils.mo
drwxr-xr-x	root/root	usr/share/locale/de/
drwxr-xr-x	root/root	usr/share/locale/de/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/de/LC_MESSAGES/libdvbv5.mo
-rw-r--r--	root/root	usr/share/locale/de/LC_MESSAGES/v4l-utils.mo
drwxr-xr-x	root/root	usr/share/locale/fr/
drwxr-xr-x	root/root	usr/share/locale/fr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fr/LC_MESSAGES/v4l-utils.mo
drwxr-xr-x	root/root	usr/share/locale/pt_BR/
drwxr-xr-x	root/root	usr/share/locale/pt_BR/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pt_BR/LC_MESSAGES/libdvbv5.mo
-rw-r--r--	root/root	usr/share/locale/pt_BR/LC_MESSAGES/v4l-utils.mo
drwxr-xr-x	root/root	usr/share/locale/uk/
drwxr-xr-x	root/root	usr/share/locale/uk/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/uk/LC_MESSAGES/libdvbv5.mo
-rw-r--r--	root/root	usr/share/locale/uk/LC_MESSAGES/v4l-utils.mo
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/cec-compliance.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cec-ctl.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cec-follower.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dvb-fe-tool.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dvb-format-convert.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dvbv5-scan.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dvbv5-zap.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ir-ctl.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ir-keytable.1.gz
-rw-r--r--	root/root	usr/share/man/man1/qvidcap.1.gz
-rw-r--r--	root/root	usr/share/man/man1/v4l2-compliance.1.gz
-rw-r--r--	root/root	usr/share/man/man1/v4l2-ctl.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/rc_keymap.5.gz
