# ts/build/packages/ica/build/extra/etc/xdg/autostart/ica_receiver_config.desktop
[Desktop Entry]
Version=1.0
Name=ica_receiver_config
Comment=Initial configuration for Citrix Receiver
TryExec=/bin/ica_receiver_config.sh
Exec=/bin/ica_receiver_config.sh
StartupNotify=false
Terminal=false
Type=Application
Icon=ica
Categories=GTK;GNOME;X-Xfce-Toplevel;