# ts/ports/core/iputils/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/arping
-rwxr-xr-x	root/root	usr/bin/clockdiff
-rwxr-xr-x	root/root	usr/bin/ping
-rwxr-xr-x	root/root	usr/bin/tracepath
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/ninfod.service
-rw-r--r--	root/root	usr/lib/systemd/system/rdisc.service
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/ninfod
-rwxr-xr-x	root/root	usr/sbin/rdisc
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/ja/
drwxr-xr-x	root/root	usr/share/locale/ja/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ja/LC_MESSAGES/iputils.mo
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/arping.8.gz
-rw-r--r--	root/root	usr/share/man/man8/clockdiff.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ninfod.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ping.8.gz
-rw-r--r--	root/root	usr/share/man/man8/rdisc.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tracepath.8.gz
