# ts/ports/core/inetutils/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/inetd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/inetd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ftp
-rwxr-xr-x	root/root	usr/bin/telnet
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/ftp.1.gz
-rw-r--r--	root/root	usr/man/man1/telnet.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/inetd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/inetd
