# ts/ports/core/inetutils/.footprint
drwxr-xr-x	root/root	bin/
-rwxr-xr-x	root/root	bin/hostname
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/inetd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/inetd
drwxr-xr-x	root/root	sbin/
-rwxr-xr-x	root/root	sbin/ifconfig
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ftp
-rwxr-xr-x	root/root	usr/bin/telnet
-rwsr-xr-x	root/root	usr/bin/traceroute
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/inetd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/ftp.1.gz
-rw-r--r--	root/root	usr/share/man/man1/hostname.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ifconfig.1.gz
-rw-r--r--	root/root	usr/share/man/man1/telnet.1.gz
-rw-r--r--	root/root	usr/share/man/man1/traceroute.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/inetd.8.gz
