# ts/build/packages/rngd/etc/systemd/system/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
Type=oneshot
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/rngd init

[Install]
WantedBy=multi-user.target

