# ts/ports/xorg-opt/accountsservice/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/org.freedesktop.Accounts.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/accountsservice-1.0/
drwxr-xr-x	root/root	usr/include/accountsservice-1.0/act/
-rw-r--r--	root/root	usr/include/accountsservice-1.0/act/act-user-enum-types.h
-rw-r--r--	root/root	usr/include/accountsservice-1.0/act/act-user-manager.h
-rw-r--r--	root/root	usr/include/accountsservice-1.0/act/act-user.h
-rw-r--r--	root/root	usr/include/accountsservice-1.0/act/act.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/girepository-1.0/
-rw-r--r--	root/root	usr/lib/girepository-1.0/AccountsService-1.0.typelib
-rwxr-xr-x	root/root	usr/lib/libaccountsservice.la
lrwxrwxrwx	root/root	usr/lib/libaccountsservice.so -> libaccountsservice.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libaccountsservice.so.0 -> libaccountsservice.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libaccountsservice.so.0.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/accountsservice.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/accounts-daemon.service
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/accounts-daemon
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/interfaces/
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.Accounts.User.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/org.freedesktop.Accounts.xml
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.Accounts.service
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/accountsservice/
drwxr-xr-x	root/root	usr/share/doc/accountsservice/spec/
-rw-r--r--	root/root	usr/share/doc/accountsservice/spec/AccountsService.html
drwxr-xr-x	root/root	usr/share/gir-1.0/
-rw-r--r--	root/root	usr/share/gir-1.0/AccountsService-1.0.gir
drwxr-xr-x	root/root	usr/share/gtk-doc/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/libaccountsservice/
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/ActUser.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/ActUserManager.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/annotation-glossary.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/api-index-full.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/ch01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/left-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/libaccountsservice.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/right-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/up-insensitive.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libaccountsservice/up.png
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/ar/
drwxr-xr-x	root/root	usr/share/locale/ar/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ar/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/bg_BG/
drwxr-xr-x	root/root	usr/share/locale/bg_BG/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/bg_BG/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/bn_IN/
drwxr-xr-x	root/root	usr/share/locale/bn_IN/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/bn_IN/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ca/
drwxr-xr-x	root/root	usr/share/locale/ca/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ca/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ca@valencia/
drwxr-xr-x	root/root	usr/share/locale/ca@valencia/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ca@valencia/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/cs/
drwxr-xr-x	root/root	usr/share/locale/cs/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/cs/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/da/
drwxr-xr-x	root/root	usr/share/locale/da/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/da/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/de/
drwxr-xr-x	root/root	usr/share/locale/de/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/de/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/el/
drwxr-xr-x	root/root	usr/share/locale/el/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/el/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/en_GB/
drwxr-xr-x	root/root	usr/share/locale/en_GB/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/en_GB/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/eo/
drwxr-xr-x	root/root	usr/share/locale/eo/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/eo/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/es/
drwxr-xr-x	root/root	usr/share/locale/es/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/es/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/fa_IR/
drwxr-xr-x	root/root	usr/share/locale/fa_IR/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fa_IR/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/fi/
drwxr-xr-x	root/root	usr/share/locale/fi/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fi/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/fr/
drwxr-xr-x	root/root	usr/share/locale/fr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fr/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ga/
drwxr-xr-x	root/root	usr/share/locale/ga/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ga/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/gl/
drwxr-xr-x	root/root	usr/share/locale/gl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/gl/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/he/
drwxr-xr-x	root/root	usr/share/locale/he/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/he/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/hi/
drwxr-xr-x	root/root	usr/share/locale/hi/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/hi/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/hu/
drwxr-xr-x	root/root	usr/share/locale/hu/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/hu/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ia/
drwxr-xr-x	root/root	usr/share/locale/ia/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ia/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/id/
drwxr-xr-x	root/root	usr/share/locale/id/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/id/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/it/
drwxr-xr-x	root/root	usr/share/locale/it/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/it/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ja/
drwxr-xr-x	root/root	usr/share/locale/ja/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ja/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ka/
drwxr-xr-x	root/root	usr/share/locale/ka/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ka/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/kk/
drwxr-xr-x	root/root	usr/share/locale/kk/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/kk/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ko/
drwxr-xr-x	root/root	usr/share/locale/ko/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ko/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/lt/
drwxr-xr-x	root/root	usr/share/locale/lt/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/lt/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/lv/
drwxr-xr-x	root/root	usr/share/locale/lv/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/lv/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/nl/
drwxr-xr-x	root/root	usr/share/locale/nl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/nl/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/pa/
drwxr-xr-x	root/root	usr/share/locale/pa/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pa/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/pl/
drwxr-xr-x	root/root	usr/share/locale/pl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pl/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/pt_BR/
drwxr-xr-x	root/root	usr/share/locale/pt_BR/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/pt_BR/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ro/
drwxr-xr-x	root/root	usr/share/locale/ro/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ro/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/ru/
drwxr-xr-x	root/root	usr/share/locale/ru/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/ru/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/sk/
drwxr-xr-x	root/root	usr/share/locale/sk/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sk/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/sl/
drwxr-xr-x	root/root	usr/share/locale/sl/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sl/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/sr/
drwxr-xr-x	root/root	usr/share/locale/sr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sr/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/sr@latin/
drwxr-xr-x	root/root	usr/share/locale/sr@latin/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sr@latin/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/sv/
drwxr-xr-x	root/root	usr/share/locale/sv/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/sv/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/te/
drwxr-xr-x	root/root	usr/share/locale/te/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/te/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/tr/
drwxr-xr-x	root/root	usr/share/locale/tr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/tr/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/uk/
drwxr-xr-x	root/root	usr/share/locale/uk/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/uk/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/vi/
drwxr-xr-x	root/root	usr/share/locale/vi/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/vi/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/wa/
drwxr-xr-x	root/root	usr/share/locale/wa/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/wa/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/zh_CN/
drwxr-xr-x	root/root	usr/share/locale/zh_CN/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/zh_CN/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/locale/zh_TW/
drwxr-xr-x	root/root	usr/share/locale/zh_TW/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/zh_TW/LC_MESSAGES/accounts-service.mo
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/org.freedesktop.accounts.policy
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/AccountsService/
drwxr-xr-x	root/root	var/lib/AccountsService/icons/
drwxr-xr-x	root/root	var/lib/AccountsService/users/
