# ts/build/packages/mjpg-streamer/build/extra/etc/systemd/system/mjpg-streamer.service
[Unit]
Description=Thistation mjpg-streamer
After=profile-setup.service pkg.service
Before=network-pre.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/mjpg-streamer init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
