# usr/ports/opt/alsa-utils/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/alsa
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/aconnect
-rwxr-xr-x	root/root	usr/bin/alsaloop
-rwxr-xr-x	root/root	usr/bin/alsamixer
-rwxr-xr-x	root/root	usr/bin/alsaucm
-rwxr-xr-x	root/root	usr/bin/amidi
-rwxr-xr-x	root/root	usr/bin/amixer
-rwxr-xr-x	root/root	usr/bin/aplay
-rwxr-xr-x	root/root	usr/bin/aplaymidi
lrwxrwxrwx	root/root	usr/bin/arecord -> aplay
-rwxr-xr-x	root/root	usr/bin/arecordmidi
-rwxr-xr-x	root/root	usr/bin/aseqdump
-rwxr-xr-x	root/root	usr/bin/aseqnet
-rwxr-xr-x	root/root	usr/bin/iecset
-rwxr-xr-x	root/root	usr/bin/speaker-test
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/aconnect.1.gz
-rw-r--r--	root/root	usr/man/man1/alsactl.1.gz
-rw-r--r--	root/root	usr/man/man1/alsaloop.1.gz
-rw-r--r--	root/root	usr/man/man1/alsamixer.1.gz
-rw-r--r--	root/root	usr/man/man1/amidi.1.gz
-rw-r--r--	root/root	usr/man/man1/amixer.1.gz
-rw-r--r--	root/root	usr/man/man1/aplay.1.gz
-rw-r--r--	root/root	usr/man/man1/aplaymidi.1.gz
lrwxrwxrwx	root/root	usr/man/man1/arecord.1.gz -> aplay.1.gz
-rw-r--r--	root/root	usr/man/man1/arecordmidi.1.gz
-rw-r--r--	root/root	usr/man/man1/aseqdump.1.gz
-rw-r--r--	root/root	usr/man/man1/aseqnet.1.gz
-rw-r--r--	root/root	usr/man/man1/iecset.1.gz
-rw-r--r--	root/root	usr/man/man1/speaker-test.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/alsaconf.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/alsaconf
-rwxr-xr-x	root/root	usr/sbin/alsactl
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/alsa/
