# ts/build/packages/nfs/etc/systemd/system/nfs.service
[Unit]
Description=ThinStation nfs
After=profile-setup.service
Before=network-pre.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/nfs init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
