# ts/ports/opt/bluez/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/bluetooth/
-rw-r--r--	root/root	etc/bluetooth/main.conf
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/bluetooth.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/advtest
-rwxr-xr-x	root/root	usr/bin/amptest
-rwxr-xr-x	root/root	usr/bin/avinfo
-rwxr-xr-x	root/root	usr/bin/avtest
-rwxr-xr-x	root/root	usr/bin/bccmd
-rwxr-xr-x	root/root	usr/bin/bcmfw
-rwxr-xr-x	root/root	usr/bin/bdaddr
-rwxr-xr-x	root/root	usr/bin/bluemoon
-rwxr-xr-x	root/root	usr/bin/bluetooth-player
-rwxr-xr-x	root/root	usr/bin/bluetoothctl
-rwxr-xr-x	root/root	usr/bin/bneptest
-rwxr-xr-x	root/root	usr/bin/btattach
-rwxr-xr-x	root/root	usr/bin/btconfig
-rwxr-xr-x	root/root	usr/bin/btgatt-client
-rwxr-xr-x	root/root	usr/bin/btgatt-server
-rwxr-xr-x	root/root	usr/bin/btinfo
-rwxr-xr-x	root/root	usr/bin/btiotest
-rwxr-xr-x	root/root	usr/bin/btmgmt
-rwxr-xr-x	root/root	usr/bin/btmon
-rwxr-xr-x	root/root	usr/bin/btproxy
-rwxr-xr-x	root/root	usr/bin/btsnoop
-rwxr-xr-x	root/root	usr/bin/check-selftest
-rwxr-xr-x	root/root	usr/bin/cltest
-rwxr-xr-x	root/root	usr/bin/create-image
-rwxr-xr-x	root/root	usr/bin/eddystone
-rwxr-xr-x	root/root	usr/bin/gatt-service
-rwxr-xr-x	root/root	usr/bin/hcieventmask
-rwxr-xr-x	root/root	usr/bin/hcisecfilter
-rwxr-xr-x	root/root	usr/bin/hex2hcd
-rwxr-xr-x	root/root	usr/bin/hid2hci
-rwxr-xr-x	root/root	usr/bin/hwdb
-rwxr-xr-x	root/root	usr/bin/ibeacon
-rwxr-xr-x	root/root	usr/bin/l2ping
-rwxr-xr-x	root/root	usr/bin/l2test
-rwxr-xr-x	root/root	usr/bin/mcaptest
-rwxr-xr-x	root/root	usr/bin/mpris-proxy
-rwxr-xr-x	root/root	usr/bin/nokfw
-rwxr-xr-x	root/root	usr/bin/obex-client-tool
-rwxr-xr-x	root/root	usr/bin/obex-server-tool
-rwxr-xr-x	root/root	usr/bin/obexctl
-rwxr-xr-x	root/root	usr/bin/oobtest
-rwxr-xr-x	root/root	usr/bin/rctest
-rwxr-xr-x	root/root	usr/bin/rtlfw
-rwxr-xr-x	root/root	usr/bin/scotest
-rwxr-xr-x	root/root	usr/bin/seq2bseq
-rwxr-xr-x	root/root	usr/bin/test-runner
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/bluetooth/
-rw-r--r--	root/root	usr/include/bluetooth/bluetooth.h
-rw-r--r--	root/root	usr/include/bluetooth/bnep.h
-rw-r--r--	root/root	usr/include/bluetooth/cmtp.h
-rw-r--r--	root/root	usr/include/bluetooth/hci.h
-rw-r--r--	root/root	usr/include/bluetooth/hci_lib.h
-rw-r--r--	root/root	usr/include/bluetooth/hidp.h
-rw-r--r--	root/root	usr/include/bluetooth/l2cap.h
-rw-r--r--	root/root	usr/include/bluetooth/rfcomm.h
-rw-r--r--	root/root	usr/include/bluetooth/sco.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp_lib.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/bluetooth/
-rwxr-xr-x	root/root	usr/lib/bluetooth/bluetoothd
-rwxr-xr-x	root/root	usr/lib/bluetooth/obexd
drwxr-xr-x	root/root	usr/lib/bluetooth/plugins/
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/sixaxis.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/sixaxis.so
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/backend/
-rwxr-xr-x	root/root	usr/lib/cups/backend/bluetooth
-rwxr-xr-x	root/root	usr/lib/libbluetooth.la
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so -> libbluetooth.so.3.19.0
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so.3 -> libbluetooth.so.3.19.0
-rwxr-xr-x	root/root	usr/lib/libbluetooth.so.3.19.0
drwxr-xr-x	root/root	usr/lib/modprobe.d/
-rw-r--r--	root/root	usr/lib/modprobe.d/bluetooth-usb.conf
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/bluez.pc
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/bluetooth.service
drwxr-xr-x	root/root	usr/lib/systemd/user/
lrwxrwxrwx	root/root	usr/lib/systemd/user/dbus-org.bluez.obex.service -> /usr/lib/systemd/user/obex.service
-rw-r--r--	root/root	usr/lib/systemd/user/obex.service
drwxr-xr-x	root/root	usr/lib/udev/
-rwxr-xr-x	root/root	usr/lib/udev/hid2hci
drwxr-xr-x	root/root	usr/lib/udev/rules.d/
-rw-r--r--	root/root	usr/lib/udev/rules.d/97-hid2hci.rules
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/services/
-rw-r--r--	root/root	usr/share/dbus-1/services/org.bluez.obex.service
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.bluez.service
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/bluez/
drwxr-xr-x	root/root	usr/share/doc/bluez/dbus-apis/
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/adapter-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/advertising-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/agent-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/assigned-numbers.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/btmon.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/btsnoop.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/device-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/gatt-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/health-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/input-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/media-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/mgmt-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/network-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/obex-agent-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/obex-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/pics-opp.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/pixit-opp.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/profile-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/pts-opp.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/sap-api.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/settings-storage.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/supported-features.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/test-coverage.txt
-rw-r--r--	root/root	usr/share/doc/bluez/dbus-apis/test-runner.txt
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/bccmd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/btattach.1.gz
-rw-r--r--	root/root	usr/share/man/man1/hid2hci.1.gz
-rw-r--r--	root/root	usr/share/man/man1/l2ping.1.gz
-rw-r--r--	root/root	usr/share/man/man1/rctest.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/bluetoothd.8.gz
