# ts/ports/core/sysklogd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/cron/
drwxr-xr-x	root/root	etc/cron/weekly/
-rwxr-xr-x	root/root	etc/cron/weekly/syslog
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/sysklogd
-rw-r--r--	root/root	etc/syslog.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/syslog.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/klogd.8.gz
-rw-r--r--	root/root	usr/man/man8/sysklogd.8.gz
-rw-r--r--	root/root	usr/man/man8/syslogd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwx------	root/root	usr/sbin/klogd
-rwxr-xr-x	root/root	usr/sbin/rotatelog
-rwx------	root/root	usr/sbin/syslogd
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
-rw-r-----	root/root	var/log/auth (EMPTY)
-rw-r-----	root/root	var/log/cron (EMPTY)
-rw-r-----	root/root	var/log/debug (EMPTY)
-rw-r-----	root/root	var/log/kernel (EMPTY)
-rw-r-----	root/root	var/log/mail (EMPTY)
-rw-r-----	root/root	var/log/messages (EMPTY)
