# ts/ports/static-source/xinetd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rw-rw-r--	root/root	etc/xinetd.conf.sample
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/xinetd.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/itox.8.gz
-rw-r--r--	root/root	usr/man/man8/xconv.pl.8.gz
-rw-r--r--	root/root	usr/man/man8/xinetd.8.gz
-rw-r--r--	root/root	usr/man/man8/xinetd.log.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/itox
-rwxr-xr-x	root/root	usr/sbin/xconv.pl
-rwxr-xr-x	root/root	usr/sbin/xinetd
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
-rw-------	root/root	var/log/xinetd.log (EMPTY)
