# ts/ports/opt/spice-vdagent/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/conf.d/
-rw-r--r--	root/root	etc/conf.d/spice-vdagentd
drwxr-xr-x	root/root	etc/xdg/
drwxr-xr-x	root/root	etc/xdg/autostart/
-rw-r--r--	root/root	etc/xdg/autostart/spice-vdagent.desktop
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/spice-vdagent
-rwxr-xr-x	root/root	usr/bin/spice-vdagentd
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/spice-vdagentd.service
-rw-r--r--	root/root	usr/lib/systemd/system/spice-vdagentd.socket
drwxr-xr-x	root/root	usr/lib/tmpfiles.d/
-rw-r--r--	root/root	usr/lib/tmpfiles.d/spice-vdagentd.conf
drwxr-xr-x	root/root	usr/lib/udev/
drwxr-xr-x	root/root	usr/lib/udev/rules.d/
-rw-r--r--	root/root	usr/lib/udev/rules.d/70-spice-vdagentd.rules
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/gdm/
drwxr-xr-x	root/root	usr/share/gdm/autostart/
drwxr-xr-x	root/root	usr/share/gdm/autostart/LoginWindow/
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/spice-vdagent.desktop
drwxr-xr-x	root/root	usr/share/gdm/greeter/
drwxr-xr-x	root/root	usr/share/gdm/greeter/autostart/
-rw-r--r--	root/root	usr/share/gdm/greeter/autostart/spice-vdagent.desktop
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/spice-vdagent.1.gz
-rw-r--r--	root/root	usr/share/man/man1/spice-vdagentd.1.gz
