# ts/ports/xorg/mesa3d/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/OpenCL/
drwxr-xr-x	root/root	etc/OpenCL/vendors/
-rw-r--r--	root/root	etc/OpenCL/vendors/mesa.icd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-r-xr-xr-x	root/root	usr/bin/mesa-overlay-control.py
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/EGL/
-rw-r--r--	root/root	usr/include/EGL/eglextchromium.h
-rw-r--r--	root/root	usr/include/EGL/eglmesaext.h
drwxr-xr-x	root/root	usr/include/GL/
drwxr-xr-x	root/root	usr/include/GL/internal/
-rw-r--r--	root/root	usr/include/GL/internal/dri_interface.h
-rw-r--r--	root/root	usr/include/GL/osmesa.h
drwxr-xr-x	root/root	usr/include/d3dadapter/
-rw-r--r--	root/root	usr/include/d3dadapter/d3dadapter9.h
-rw-r--r--	root/root	usr/include/d3dadapter/drm.h
-rw-r--r--	root/root	usr/include/d3dadapter/present.h
-rw-r--r--	root/root	usr/include/gbm.h
-rw-r--r--	root/root	usr/include/xa_composite.h
-rw-r--r--	root/root	usr/include/xa_context.h
-rw-r--r--	root/root	usr/include/xa_tracker.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/bellagio/
-rwxr-xr-x	root/root	usr/lib/bellagio/libomx_mesa.so
drwxr-xr-x	root/root	usr/lib/d3d/
lrwxrwxrwx	root/root	usr/lib/d3d/d3dadapter9.so -> d3dadapter9.so.1
lrwxrwxrwx	root/root	usr/lib/d3d/d3dadapter9.so.1 -> d3dadapter9.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/d3d/d3dadapter9.so.1.0.0
drwxr-xr-x	root/root	usr/lib/dri/
-rwxr-xr-x	root/root	usr/lib/dri/i830_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/i915_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/i965_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/iris_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/kms_swrast_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/nouveau_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/nouveau_drv_video.so
-rwxr-xr-x	root/root	usr/lib/dri/nouveau_vieux_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/r200_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/r300_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/r600_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/r600_drv_video.so
-rwxr-xr-x	root/root	usr/lib/dri/radeon_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/radeonsi_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/radeonsi_drv_video.so
-rwxr-xr-x	root/root	usr/lib/dri/swrast_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/virtio_gpu_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/vmwgfx_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/zink_dri.so
drwxr-xr-x	root/root	usr/lib/gallium-pipe/
-rwxr-xr-x	root/root	usr/lib/gallium-pipe/pipe_iris.so
-rwxr-xr-x	root/root	usr/lib/gallium-pipe/pipe_nouveau.so
-rwxr-xr-x	root/root	usr/lib/gallium-pipe/pipe_r300.so
-rwxr-xr-x	root/root	usr/lib/gallium-pipe/pipe_r600.so
-rwxr-xr-x	root/root	usr/lib/gallium-pipe/pipe_radeonsi.so
-rwxr-xr-x	root/root	usr/lib/gallium-pipe/pipe_swrast.so
-rwxr-xr-x	root/root	usr/lib/gallium-pipe/pipe_vmwgfx.so
lrwxrwxrwx	root/root	usr/lib/libEGL_mesa.so -> libEGL_mesa.so.0
lrwxrwxrwx	root/root	usr/lib/libEGL_mesa.so.0 -> libEGL_mesa.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libEGL_mesa.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libGLX_mesa.so -> libGLX_mesa.so.0
lrwxrwxrwx	root/root	usr/lib/libGLX_mesa.so.0 -> libGLX_mesa.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libGLX_mesa.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libMesaOpenCL.so -> libMesaOpenCL.so.1
lrwxrwxrwx	root/root	usr/lib/libMesaOpenCL.so.1 -> libMesaOpenCL.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libMesaOpenCL.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libOSMesa.so -> libOSMesa.so.8
lrwxrwxrwx	root/root	usr/lib/libOSMesa.so.8 -> libOSMesa.so.8.0.0
-rwxr-xr-x	root/root	usr/lib/libOSMesa.so.8.0.0
-rwxr-xr-x	root/root	usr/lib/libVkLayer_INTEL_nullhw.so
-rwxr-xr-x	root/root	usr/lib/libVkLayer_MESA_device_select.so
-rwxr-xr-x	root/root	usr/lib/libVkLayer_MESA_overlay.so
lrwxrwxrwx	root/root	usr/lib/libXvMCnouveau.so -> libXvMCnouveau.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libXvMCnouveau.so.1 -> libXvMCnouveau.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libXvMCnouveau.so.1.0 -> libXvMCnouveau.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libXvMCnouveau.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libXvMCr600.so -> libXvMCr600.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libXvMCr600.so.1 -> libXvMCr600.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libXvMCr600.so.1.0 -> libXvMCr600.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libXvMCr600.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libgbm.so -> libgbm.so.1
lrwxrwxrwx	root/root	usr/lib/libgbm.so.1 -> libgbm.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libgbm.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libglapi.so -> libglapi.so.0
lrwxrwxrwx	root/root	usr/lib/libglapi.so.0 -> libglapi.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libglapi.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libswrAVX.so -> libswrAVX.so.0
lrwxrwxrwx	root/root	usr/lib/libswrAVX.so.0 -> libswrAVX.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libswrAVX.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libswrAVX2.so -> libswrAVX2.so.0
lrwxrwxrwx	root/root	usr/lib/libswrAVX2.so.0 -> libswrAVX2.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libswrAVX2.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libvulkan_intel.so
-rwxr-xr-x	root/root	usr/lib/libvulkan_lvp.so
-rwxr-xr-x	root/root	usr/lib/libvulkan_radeon.so
lrwxrwxrwx	root/root	usr/lib/libxatracker.so -> libxatracker.so.2
lrwxrwxrwx	root/root	usr/lib/libxatracker.so.2 -> libxatracker.so.2.5.0
-rwxr-xr-x	root/root	usr/lib/libxatracker.so.2.5.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/d3d.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/dri.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/gbm.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/osmesa.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/xatracker.pc
drwxr-xr-x	root/root	usr/lib/vdpau/
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_nouveau.so -> libvdpau_nouveau.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_nouveau.so.1 -> libvdpau_nouveau.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_nouveau.so.1.0 -> libvdpau_nouveau.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/vdpau/libvdpau_nouveau.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_r300.so -> libvdpau_r300.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_r300.so.1 -> libvdpau_r300.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_r300.so.1.0 -> libvdpau_r300.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/vdpau/libvdpau_r300.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_r600.so -> libvdpau_r600.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_r600.so.1 -> libvdpau_r600.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_r600.so.1.0 -> libvdpau_r600.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/vdpau/libvdpau_r600.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_radeonsi.so -> libvdpau_radeonsi.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_radeonsi.so.1 -> libvdpau_radeonsi.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/vdpau/libvdpau_radeonsi.so.1.0 -> libvdpau_radeonsi.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/vdpau/libvdpau_radeonsi.so.1.0.0
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/drirc.d/
-rw-r--r--	root/root	usr/share/drirc.d/00-mesa-defaults.conf
drwxr-xr-x	root/root	usr/share/glvnd/
drwxr-xr-x	root/root	usr/share/glvnd/egl_vendor.d/
-rw-r--r--	root/root	usr/share/glvnd/egl_vendor.d/50_mesa.json
drwxr-xr-x	root/root	usr/share/vulkan/
drwxr-xr-x	root/root	usr/share/vulkan/explicit_layer.d/
-rw-r--r--	root/root	usr/share/vulkan/explicit_layer.d/VkLayer_INTEL_nullhw.json
-rw-r--r--	root/root	usr/share/vulkan/explicit_layer.d/VkLayer_MESA_overlay.json
drwxr-xr-x	root/root	usr/share/vulkan/icd.d/
-rw-r--r--	root/root	usr/share/vulkan/icd.d/intel_icd.x86_64.json
-rw-r--r--	root/root	usr/share/vulkan/icd.d/lvp_icd.x86_64.json
-rw-r--r--	root/root	usr/share/vulkan/icd.d/radeon_icd.x86_64.json
drwxr-xr-x	root/root	usr/share/vulkan/implicit_layer.d/
-rw-r--r--	root/root	usr/share/vulkan/implicit_layer.d/VkLayer_MESA_device_select.json
