# ts/6.0/packages/networkmanager/build/extra/etc/udev/rules.d/99-tsnet.rules
SUBSYSTEM!="net", GOTO="net_end"
ACTION=="add", RUN+="/etc/udev/scripts/pre-net.sh"
LABEL="net_end"
