# ts/ports/opt/mysql/.footprint
drwxr-xr-x	root/root	etc/
-rw-------	root/root	etc/my.cnf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/mysqld
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/innochecksum
-rwxr-xr-x	root/root	usr/bin/lz4_decompress
-rwxr-xr-x	root/root	usr/bin/my_print_defaults
-rwxr-xr-x	root/root	usr/bin/myisam_ftdump
-rwxr-xr-x	root/root	usr/bin/myisamchk
-rwxr-xr-x	root/root	usr/bin/myisamlog
-rwxr-xr-x	root/root	usr/bin/myisampack
-rwxr-xr-x	root/root	usr/bin/mysql
-rwxr-xr-x	root/root	usr/bin/mysql_client_test
-rwxr-xr-x	root/root	usr/bin/mysql_client_test_embedded
-rwxr-xr-x	root/root	usr/bin/mysql_config
-rwxr-xr-x	root/root	usr/bin/mysql_config_editor
-rwxr-xr-x	root/root	usr/bin/mysql_embedded
-rwxr-xr-x	root/root	usr/bin/mysql_install_db
-rwxr-xr-x	root/root	usr/bin/mysql_plugin
-rwxr-xr-x	root/root	usr/bin/mysql_secure_installation
-rwxr-xr-x	root/root	usr/bin/mysql_ssl_rsa_setup
-rwxr-xr-x	root/root	usr/bin/mysql_tzinfo_to_sql
-rwxr-xr-x	root/root	usr/bin/mysql_upgrade
-rwxr-xr-x	root/root	usr/bin/mysqladmin
-rwxr-xr-x	root/root	usr/bin/mysqlbinlog
-rwxr-xr-x	root/root	usr/bin/mysqlcheck
-rwxr-xr-x	root/root	usr/bin/mysqld_multi
-rwxr-xr-x	root/root	usr/bin/mysqld_safe
-rwxr-xr-x	root/root	usr/bin/mysqldump
-rwxr-xr-x	root/root	usr/bin/mysqldumpslow
-rwxr-xr-x	root/root	usr/bin/mysqlimport
-rwxr-xr-x	root/root	usr/bin/mysqlpump
-rwxr-xr-x	root/root	usr/bin/mysqlshow
-rwxr-xr-x	root/root	usr/bin/mysqlslap
-rwxr-xr-x	root/root	usr/bin/mysqltest
-rwxr-xr-x	root/root	usr/bin/mysqltest_embedded
-rwxr-xr-x	root/root	usr/bin/mysqlxtest
-rwxr-xr-x	root/root	usr/bin/perror
-rwxr-xr-x	root/root	usr/bin/replace
-rwxr-xr-x	root/root	usr/bin/resolve_stack_dump
-rwxr-xr-x	root/root	usr/bin/resolveip
-rwxr-xr-x	root/root	usr/bin/zlib_decompress
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/mysql/
-rw-r--r--	root/root	usr/include/mysql/big_endian.h
-rw-r--r--	root/root	usr/include/mysql/binary_log_types.h
-rw-r--r--	root/root	usr/include/mysql/byte_order_generic.h
-rw-r--r--	root/root	usr/include/mysql/byte_order_generic_x86.h
-rw-r--r--	root/root	usr/include/mysql/decimal.h
-rw-r--r--	root/root	usr/include/mysql/errmsg.h
-rw-r--r--	root/root	usr/include/mysql/keycache.h
-rw-r--r--	root/root	usr/include/mysql/little_endian.h
-rw-r--r--	root/root	usr/include/mysql/m_ctype.h
-rw-r--r--	root/root	usr/include/mysql/m_string.h
-rw-r--r--	root/root	usr/include/mysql/my_alloc.h
-rw-r--r--	root/root	usr/include/mysql/my_byteorder.h
-rw-r--r--	root/root	usr/include/mysql/my_command.h
-rw-r--r--	root/root	usr/include/mysql/my_compiler.h
-rw-r--r--	root/root	usr/include/mysql/my_config.h
-rw-r--r--	root/root	usr/include/mysql/my_dbug.h
-rw-r--r--	root/root	usr/include/mysql/my_dir.h
-rw-r--r--	root/root	usr/include/mysql/my_getopt.h
-rw-r--r--	root/root	usr/include/mysql/my_global.h
-rw-r--r--	root/root	usr/include/mysql/my_list.h
-rw-r--r--	root/root	usr/include/mysql/my_sys.h
-rw-r--r--	root/root	usr/include/mysql/my_thread.h
-rw-r--r--	root/root	usr/include/mysql/my_thread_local.h
-rw-r--r--	root/root	usr/include/mysql/my_xml.h
-rw-r--r--	root/root	usr/include/mysql/mysql.h
drwxr-xr-x	root/root	usr/include/mysql/mysql/
-rw-r--r--	root/root	usr/include/mysql/mysql/client_authentication.h
-rw-r--r--	root/root	usr/include/mysql/mysql/client_plugin.h
-rw-r--r--	root/root	usr/include/mysql/mysql/client_plugin.h.pp
-rw-r--r--	root/root	usr/include/mysql/mysql/com_data.h
-rw-r--r--	root/root	usr/include/mysql/mysql/get_password.h
-rw-r--r--	root/root	usr/include/mysql/mysql/group_replication_priv.h
-rw-r--r--	root/root	usr/include/mysql/mysql/innodb_priv.h
-rw-r--r--	root/root	usr/include/mysql/mysql/mysql_lex_string.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_audit.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_audit.h.pp
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_auth.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_auth.h.pp
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_auth_common.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_ftparser.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_ftparser.h.pp
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_group_replication.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_keyring.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_keyring.h.pp
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_trace.h
-rw-r--r--	root/root	usr/include/mysql/mysql/plugin_validate_password.h
drwxr-xr-x	root/root	usr/include/mysql/mysql/psi/
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_file.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_idle.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_mdl.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_memory.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_ps.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_socket.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_sp.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_stage.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_statement.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_table.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_thread.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/mysql_transaction.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/psi.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/psi_base.h
-rw-r--r--	root/root	usr/include/mysql/mysql/psi/psi_memory.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_command.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_locking.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_my_plugin_log.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_my_snprintf.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_mysql_alloc.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_mysql_keyring.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_mysql_password_policy.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_mysql_string.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_parser.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_rpl_transaction_ctx.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_rpl_transaction_write_set.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_rules_table.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_security_context.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_srv_session.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_srv_session_info.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_ssl_wrapper.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_thd_alloc.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_thd_engine_lock.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_thd_wait.h
-rw-r--r--	root/root	usr/include/mysql/mysql/service_thread_scheduler.h
-rw-r--r--	root/root	usr/include/mysql/mysql/services.h
-rw-r--r--	root/root	usr/include/mysql/mysql/services.h.pp
-rw-r--r--	root/root	usr/include/mysql/mysql/thread_pool_priv.h
-rw-r--r--	root/root	usr/include/mysql/mysql/thread_type.h
-rw-r--r--	root/root	usr/include/mysql/mysql_com.h
-rw-r--r--	root/root	usr/include/mysql/mysql_com_server.h
-rw-r--r--	root/root	usr/include/mysql/mysql_embed.h
-rw-r--r--	root/root	usr/include/mysql/mysql_time.h
-rw-r--r--	root/root	usr/include/mysql/mysql_version.h
-rw-r--r--	root/root	usr/include/mysql/mysqld_ername.h
-rw-r--r--	root/root	usr/include/mysql/mysqld_error.h
-rw-r--r--	root/root	usr/include/mysql/mysqlx_ername.h
-rw-r--r--	root/root	usr/include/mysql/mysqlx_error.h
-rw-r--r--	root/root	usr/include/mysql/mysqlx_version.h
-rw-r--r--	root/root	usr/include/mysql/plugin.h
-rw-r--r--	root/root	usr/include/mysql/plugin_audit.h
-rw-r--r--	root/root	usr/include/mysql/plugin_ftparser.h
-rw-r--r--	root/root	usr/include/mysql/plugin_group_replication.h
-rw-r--r--	root/root	usr/include/mysql/plugin_keyring.h
-rw-r--r--	root/root	usr/include/mysql/plugin_validate_password.h
-rw-r--r--	root/root	usr/include/mysql/sql_common.h
-rw-r--r--	root/root	usr/include/mysql/sql_state.h
-rw-r--r--	root/root	usr/include/mysql/sslopt-case.h
-rw-r--r--	root/root	usr/include/mysql/sslopt-longopts.h
-rw-r--r--	root/root	usr/include/mysql/sslopt-vars.h
-rw-r--r--	root/root	usr/include/mysql/thr_cond.h
-rw-r--r--	root/root	usr/include/mysql/thr_mutex.h
-rw-r--r--	root/root	usr/include/mysql/thr_rwlock.h
-rw-r--r--	root/root	usr/include/mysql/typelib.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libmysqlclient.a
lrwxrwxrwx	root/root	usr/lib/libmysqlclient.so -> libmysqlclient.so.20
lrwxrwxrwx	root/root	usr/lib/libmysqlclient.so.20 -> libmysqlclient.so.20.3.14
-rwxr-xr-x	root/root	usr/lib/libmysqlclient.so.20.3.14
-rw-r--r--	root/root	usr/lib/libmysqld.a
-rw-r--r--	root/root	usr/lib/libmysqlservices.a
drwxr-xr-x	root/root	usr/lib/mysql/
drwxr-xr-x	root/root	usr/lib/mysql/plugin/
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/adt_null.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/auth.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/auth_socket.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/auth_test_plugin.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/authentication_ldap_sasl_client.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/connection_control.so
-rw-r--r--	root/root	usr/lib/mysql/plugin/daemon_example.ini
drwxr-xr-x	root/root	usr/lib/mysql/plugin/debug/
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/group_replication.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/ha_example.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/keyring_file.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/keyring_udf.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libdaemon_example.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_framework.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_services.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_services_threaded.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_session_detach.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_session_in_thd.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_session_info.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_2_sessions.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_all_col_types.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_cmds_1.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_commit.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_complex.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_errors.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_lock.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_processlist.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_replication.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_shutdown.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_sqlmode.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_stored_procedures_functions.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_sql_views_triggers.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_x_sessions_deinit.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/libtest_x_sessions_init.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/locking_service.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/mypluglib.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/mysql_no_login.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/mysqlx.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/qa_auth_client.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/qa_auth_interface.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/qa_auth_server.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/replication_observers_example_plugin.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/rewrite_example.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/rewriter.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/semisync_master.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/semisync_slave.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/test_security_context.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/test_udf_services.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/udf_example.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/validate_password.so
-rwxr-xr-x	root/root	usr/lib/mysql/plugin/version_token.so
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/mysqlclient.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/mysqld
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/aclocal/
-rw-r--r--	root/root	usr/share/aclocal/mysql.m4
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/comp_err.1.gz
-rw-r--r--	root/root	usr/share/man/man1/innochecksum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/lz4_decompress.1.gz
-rw-r--r--	root/root	usr/share/man/man1/my_print_defaults.1.gz
-rw-r--r--	root/root	usr/share/man/man1/myisam_ftdump.1.gz
-rw-r--r--	root/root	usr/share/man/man1/myisamchk.1.gz
-rw-r--r--	root/root	usr/share/man/man1/myisamlog.1.gz
-rw-r--r--	root/root	usr/share/man/man1/myisampack.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql.server.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_config.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_config_editor.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_install_db.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_plugin.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_secure_installation.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_ssl_rsa_setup.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_tzinfo_to_sql.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysql_upgrade.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqladmin.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqlbinlog.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqlcheck.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqld_multi.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqld_safe.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqldump.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqldumpslow.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqlimport.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqlman.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqlpump.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqlshow.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mysqlslap.1.gz
-rw-r--r--	root/root	usr/share/man/man1/perror.1.gz
-rw-r--r--	root/root	usr/share/man/man1/replace.1.gz
-rw-r--r--	root/root	usr/share/man/man1/resolve_stack_dump.1.gz
-rw-r--r--	root/root	usr/share/man/man1/resolveip.1.gz
-rw-r--r--	root/root	usr/share/man/man1/zlib_decompress.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/mysqld.8.gz
drwxr-xr-x	root/root	usr/share/mysql/
drwxr-xr-x	root/root	usr/share/mysql/charsets/
-rw-r--r--	root/root	usr/share/mysql/charsets/Index.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/armscii8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/ascii.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1250.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1251.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1256.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp1257.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp850.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp852.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/cp866.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/dec8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/geostd8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/greek.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/hebrew.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/hp8.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/keybcs2.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/koi8r.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/koi8u.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin1.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin2.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin5.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/latin7.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/macce.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/macroman.xml
-rw-r--r--	root/root	usr/share/mysql/charsets/swe7.xml
-rw-r--r--	root/root	usr/share/mysql/dictionary.txt
drwxr-xr-x	root/root	usr/share/mysql/english/
-rw-r--r--	root/root	usr/share/mysql/english/errmsg.sys
-rw-r--r--	root/root	usr/share/mysql/errmsg-utf8.txt
-rw-r--r--	root/root	usr/share/mysql/fill_help_tables.sql
-rw-r--r--	root/root	usr/share/mysql/innodb_memcached_config.sql
-rw-r--r--	root/root	usr/share/mysql/install_rewriter.sql
-rw-r--r--	root/root	usr/share/mysql/magic
-rwxr-xr-x	root/root	usr/share/mysql/mysql-log-rotate
-rwxr-xr-x	root/root	usr/share/mysql/mysql.server
-rw-r--r--	root/root	usr/share/mysql/mysql_security_commands.sql
-rw-r--r--	root/root	usr/share/mysql/mysql_sys_schema.sql
-rw-r--r--	root/root	usr/share/mysql/mysql_system_tables.sql
-rw-r--r--	root/root	usr/share/mysql/mysql_system_tables_data.sql
-rw-r--r--	root/root	usr/share/mysql/mysql_test_data_timezone.sql
-rwxr-xr-x	root/root	usr/share/mysql/mysqld_multi.server
-rw-r--r--	root/root	usr/share/mysql/uninstall_rewriter.sql
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwx------	mysql/mysql	var/lib/mysql-files/
drwx------	mysql/mysql	var/lib/mysql/
drwxr-xr-x	root/root	var/log/
drwx------	mysql/mysql	var/log/mysql/
drwxr-xr-x	root/root	var/run/
drwx------	mysql/mysql	var/run/mysql/
