# ts/ports/core/rc/rc.local
#!/bin/bash
#
# /etc/rc.local: local multi-user startup script
#

# End of file
