# ts/build/packages/systemd/build/extra/etc/systemd/user/homedir.service
[Unit]
Description=Create User Home Directory

[Service]
ExecStart=/sbin/session open

[Install]
WantedBy=default.target
