# ts/5.1/packages/volumeicon/etc/xdg/autostart/volumeicon.desktop
[Desktop Entry]
Name=Volumeicon
Comment=Simple volumeicon
Icon=volumeicon
Exec=runonce volumeicon
Terminal=false
Type=Application
NoDisplay=true
NotShowIn=KDE;
X-GNOME-Bugzilla-Bugzilla=GNOME
X-GNOME-Bugzilla-Product=volumeicon
X-GNOME-Bugzilla-Component=general
X-GNOME-Autostart-enabled=true
