# ts/build/packages/gnome-core/build/extra/etc/systemd/system/gnome-core.service
[Unit]
Description=Thinstation Gnome Core
Conflicts=shutdown.target
After=session.service pkg.service
Before=shutdown.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/gnome-core
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
