# ts/build/packages/set-resolution/etc/xdg/autostart/set-resolution.desktop
[Desktop Entry]
Version=1.0
Name=set-resolution
Comment=Initial configuration for set-resoultion package to extend over multiple monitors
TryExec=/bin/set-resolution
Exec=/bin/set-resolution init
StartupNotify=false
Terminal=false
Type=Application
Categories=GTK;GNOME;X-Xfce-Toplevel;
