# ts/build/packages/set-resolution/etc/xdg/autostart/set-monitor.desktop
[Desktop Entry]
Version=1.0
Name=set-monitor
Comment=Initial configuration for set-resoultion package to extend over multiple monitors
TryExec=/bin/set-monitor
Exec=/bin/set-monitor
StartupNotify=false
Terminal=false
Type=Application
Icon=ica
Categories=GTK;GNOME;X-Xfce-Toplevel;
