# usr/ports/opt/stunnel/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/stunnel
drwxr-xr-x	root/root	etc/ssl/
drwxr-xr-x	root/root	etc/ssl/certs/
-rw-r--r--	root/root	etc/ssl/certs/stunnel.crt (EMPTY)
drwxr-xr-x	root/root	etc/ssl/keys/
-rw-------	root/root	etc/ssl/keys/stunnel.key (EMPTY)
-rw-r--r--	root/root	etc/stunnel.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/stunnel.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/stunnel
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/run/
drwxrwx---	stunnel/stunnel	var/run/stunnel/
