# ts/build/packages/cups/build/extra/lib/systemd/system/cups.init.service
[Unit]
Description=Configure Thinstation Xorg Settings
DefaultDependencies=no
Conflicts=shutdown.target
After=profile_setup.service net.profle_setup.service org.cups.cupsd.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/lib/thinstation/system/cups-init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
