# ts/build/packages/devstation/etc/systemd/system/persistence-backup.service
[Unit]
Description=Backup critical files on shutdown
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target umount.target

[Service]
Type=oneshot
ExecStart=/etc/init.d/persistent-files backup
RemainAfterExit=true

[Install]
WantedBy=halt.target reboot.target shutdown.target
