# ts/build/packages/devstation/etc/systemd/system/persistence-backup.timer
[Unit]
Description=Timer to periodically back up critical system files

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
