# ts/build/packages/devstation/etc/systemd/system/persistence-restore.service
[Unit]
Description=Check and link system critical files
StartAfter=fastboot.service

[Service]
Type=oneshot
ExecStart=/etc/init.d/persistent-files restore
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
