# ts/build/packages/hdupdate/etc/systemd/system/hdupdate.service
[Unit]
Description=ThinStation hdupdate
After=network-files.service
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/etc/init.d/hdupdate init
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
