# ts/ports/kernel-modules/usb-redirector/.footprint
drwxr-xr-x	root/root	bin/
-rwxr-xr-x	root/root	bin/usbclnt
-rwxr-xr-x	root/root	bin/usbsrv
-rwxr--r--	root/root	bin/usbsrvd
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/init.d/
-rwxr-xr-x	root/root	etc/init.d/rc.usbsrvd
-rwxr-xr-x	root/root	etc/init.d/rc.usbsrvd-excludes
drwxr-xr-x	root/root	etc/rc5.d/
lrwxrwxrwx	root/root	etc/rc5.d/S10rc.usbsrvd -> /etc/init.d/rc.usbsrvd
lrwxrwxrwx	root/root	etc/rc5.d/S90rc.usbsrvd-excludes -> /etc/init.d/rc.usbsrvd-excludes
-rw-r--r--	root/root	etc/usbsrvd.conf
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/modules/
drwxr-xr-x	root/root	lib/modules/3.14.24TS/
drwxr-xr-x	root/root	lib/modules/3.14.24TS/kernel/
-rw-r--r--	root/root	lib/modules/3.14.24TS/kernel/tusbd.ko
drwxr-xr-x	root/root	lib/modules/3.14.24TS_64/
drwxr-xr-x	root/root	lib/modules/3.14.24TS_64/kernel/
-rw-r--r--	root/root	lib/modules/3.14.24TS_64/kernel/tusbd.ko
drwxr-xr-x	root/root	lib/modules/3.14.24TS_SMP/
drwxr-xr-x	root/root	lib/modules/3.14.24TS_SMP/kernel/
-rw-r--r--	root/root	lib/modules/3.14.24TS_SMP/kernel/tusbd.ko
