# ts/build/packages/base/lib/systemd/system/session.service
[Unit]
Description=Thistation Session Initialize
After=profile_setup.service
Before=shutdown.target
ConditionPathIsReadWrite=/etc

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/thinstation.env
ExecStart=/lib/thinstation/system/session
SyslogIdentifier=thinstation

[Install]
WantedBy=multi-user.target
